# bareos-common/bareos-fd.service
# systemd service file for Bareos File Deaemon
#
# we do not take the service file provided by bareos package as it needs too
# much customization

[Unit]
Description=Bareos File Deaemon
Requires=network.target
After=network.target remote-fs.target syslog.target

[Service]
Type=forking
Restart=always
ExecStart=/usr/bin/bareos-fd -c /etc/bareos/bareos-fd.conf
IOSchedulingClass=idle
PIDFile=/run/bareos-fd.9102.pid

[Install]
WantedBy=multi-user.target
