# uhr/uhr@.service
[Unit]
Description=Full-screen analog clock on %i
After=dev-tty11.device systemd-user-sessions.service

[Service]
Type=simple
ExecStart=/usr/bin/uhr
StandardOutput=tty
TTYPath=/dev/%i

[Install]
WantedBy=multi-user.target
