# qubes-vm-xen/proc-xen.automount
[Unit]
Description=Automount /proc/xen files
ConditionPathExists=/proc/xen

[Automount]
Where=/proc/xen
