# qubes-vm-xen/xenconsoled.service
[Unit]
Description=Xenconsoled - handles logging from guest consoles and hypervisor
After=xenstored.service

[Service]
Type=simple
PIDFile=/var/run/xenconsoled.pid
ExecStart=/usr/sbin/xenconsoled

[Install]
WantedBy=multi-user.target
