# qubes-vm-xen/xenstored.service
[Unit]
Description=Xenstored - daemon managing xenstore file system
Before=libvirtd.service libvirt-guests.service proc-xen.mount
After=dbus.service
RefuseManualStop=true

[Service]
Type=forking
PIDFile=/var/run/xenstored.pid
ExecStart=/usr/sbin/xenstored --pid-file /var/run/xenstored.pid $XENSTORED_ARGS

[Install]
WantedBy=multi-user.target
