# kaoz/kaoz.service
[Unit]
Description=An IRC notifier

[Service]
ExecStart=/usr/bin/kaoz

[Install]
WantedBy=multi-user.target
