# gearman/gearman.service
[Unit]
Description=Gearman server

[Service]
User=gearman
ExecStart=/usr/bin/gearmand
Restart=always

[Install]
WantedBy=multi-user.target
