# krb5-appl-git/krb5-eklogind@.service
[Unit]
Description=Encrypted Kerberized Remote Login service

[Service]
ExecStart=/usr/bin/krb5-klogind -c -e
StandardInput=socket
