# krb5-appl-git/krb5-klogind@.service
[Unit]
Description=Kerberized Remote Login service

[Service]
ExecStart=/usr/bin/krb5-klogind -c
StandardInput=socket
