# krb5-appl-git/krb5-eklogind.socket
[Unit]
Description=Encrypted Kerberized Remote Login socket

[Socket]
ListenStream=2105
Accept=yes

[Install]
WantedBy=sockets.target
