# krb5-appl-git/krb5-klogind.socket
[Unit]
Description=Kerberized Remote Login socket

[Socket]
ListenStream=543
Accept=yes

[Install]
WantedBy=sockets.target
