# krb5-appl-git/krb5-kshd@.service
[Unit]
Description=Kerberized Remote Shell server

[Service]
ExecStart=/usr/bin/krb5-kshd -5c
StandardInput=socket
