# krb5-appl-git/krb5-kshd.socket
[Unit]
Description=Kerberized Remote Shell socket

[Socket]
ListenStream=544
Accept=yes

[Install]
WantedBy=sockets.target
