# krb5-appl-git/krb5-telnetd@.service
[Unit]
Description=Kerberized Telnet server

[Service]
ExecStart=/usr/bin/krb5-telnetd -h
StandardInput=socket
