# privoxy-systemd/privoxy.service
[Unit]
Description=Privacy enhancing HTTP Proxy
 
[Service]
ExecStart=/usr/sbin/privoxy --no-daemon /etc/privoxy/config
StandardOutput=syslog
StandardError=syslog

[Install]
WantedBy=multi-user.target
