# glusterfs-git/glusterd.service
# Glusterd service file. Taken from the fedora packages provided by the gluster project
[Unit]
Description=GlusterFS a clustered file-system server
After=network.target rpcbind.service

[Service]
Type=forking
PIDFile=/run/glusterd.pid
LimitNOFILE=65536
ExecStart=/usr/bin/glusterd -p /run/glusterd.pid

[Install]
WantedBy=multi-user.target
