# yate-svn/yate.service
[Unit]
Description=Yet Another Telephony Engine
Documentation=man:yate(8)
Requires=network.target

After=network.target

[Service]
Type=forking
ExecStart=/usr/bin/yate -d
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=always

[Install]
WantedBy=multi-user.target

