# sickbeard-git/sickbeard.service
[Unit]
Description=Sick Beard Daemon

[Service]
Type=forking
User=sickbeard
ExecStart=/usr/bin/sickbeard --daemon --config /opt/sickbeard/config.ini --datadir /opt/sickbeard --pidfile /opt/sickbeard/sickbeard/sickbeard.pid

[Install]
WantedBy=multi-user.target
