# ajenti/ajenti.service
[Unit]
Description=Ajenti - an easy server administration frontend

[Service]
Type=simple
ExecStart=/usr/bin/ajenti-panel -d
KillMode=process
PIDFile=/run/ajenti-ipc.sock

[Install]
WantedBy=multi-user.target
