# xbmc-addon-steam-launcher/steam-switcher.service
[Unit]
Description=stops xbmc and launches steam-standalone, relaunches xbmc on close
After=network.target
BindsTo=steam-standalone.service
Conflicts=xbmc.service

[Service]
RemainAfterExit=yes
Type=oneshot
ExecStart=/usr/bin/sleep 1
ExecStopPost=/usr/bin/systemctl start xbmc.service
