# xbmc-addon-steam-launcher/steam-switcher.service
[Unit]
Description=switch from xbmc-standalone to steam-standalone
After=network.target
BindsTo=steam-standalone.service
Conflicts=xbmc.service

[Service]
RemainAfterExit=yes
Type=oneshot
ExecStart=/usr/bin/sleep 2
ExecStopPost=/usr/bin/systemctl start xbmc.service
