# dnssec-trigger/dnssec-triggerd.service
[Unit]
Description=Reconfigure local DNS(SEC) resolver on network change
After=syslog.target network.target
After=dnssec-triggerd-keygen.service
Wants=dnssec-triggerd-keygen.service
After=unbound.service
Wants=unbound.service

[Service]
Type=simple
ExecStart=/usr/bin/dnssec-triggerd -d -c /etc/dnssec-trigger/dnssec-trigger.conf
ExecStartPost=/etc/NetworkManager/dispatcher.d/01-dnssec-trigger-hook
Restart=always
RestartSec=0

[Install]
WantedBy=multi-user.target
