# sandfox-extras/sandfox-unmount-poweroff.service
[Unit]
Description=Close all Sandfox sandboxes

[Service]
Type=oneshot
ExecStart=/usr/bin/printf 'Closing Sandfox sandboxes'
ExecStart=/usr/bin/sandfox --user todd --closeall
TimeoutSec=5

[Install]
WantedBy=shutdown.target reboot.target halt.target poweroff.target
