# sandfox-extras/sandfox-rmdir-boot.service
[Unit]
Description=Delete previous Sandfox sandboxes (directories)

[Service]
Type=oneshot
ExecStart=/usr/bin/printf 'Deleting previous Sandfox sandboxes'
ExecStart=/usr/bin/sh -c '/usr/bin/rm -rf /mnt/sandfox/firefox*'

[Install]
WantedBy=multi-user.target
