# irggu-server-git/irggu-server.service
[Unit]
Description=IrGGu-Server
After=network.target

[Service]
Type=forking
PIDFile=/run/irggu/irggu.pid
User=irggu
Group=irggu
ExecStart=/usr/bin/irggu-server
ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/bin/kill -TERM $MAINPID

[Install]
WantedBy=multi-user.target