# scm-manager/scm-manager.service
[Unit]
Description=SCM manager

[Service]
Type=forking
ExecStart=/usr/sbin/scm-server start
ExecStop=/usr/sbin/scm-server stop

[Install]
WantedBy=multi-user.target

