# bundy/bundy.service
[Unit]
Description=bundy
After=syslog.target network.target

[Service]
ExecStart=/usr/bin/bundy

[Install]
WantedBy=multi-user.target
