# teamcity-agent/teamcity-agent.service
[Unit]
Description=TeamCity agent
After=network.target

[Service]
EnvironmentFile=/etc/conf.d/teamcity-agent
ExecStart=/opt/teamcity-agent/bin/agent.sh start
ExecStop=/opt/teamcity-agent/bin/agent.sh stop
Type=forking
Restart=no
TimeoutSec=30
RemainAfterExit=yes
User=teamcity-agent

[Install]
WantedBy=multi-user.target

