# xllmnrd/xllmnrd.service
[Unit]
Description=LLMNR Name Resolution Responder

[Service]
Type=forking
ExecStart=/usr/bin/xllmnrd
User=daemon

[Install]
WantedBy=multi-user.target
