# mopidy/mopidy@.service
[Unit]
Description=Mopidy music server
After=network.target

[Service]
Type=simple
User=%I
PAMName=system-local-login
ExecStart=/usr/bin/mopidy

[Install]
WantedBy=multi-user.target
