# mopidy/mopidy.service
[Unit]
Description=mopidy
After=network.target

[Service]
User=mopidy
ExecStart=/usr/bin/mopidy --config /etc/mopidy/extensions.d:/etc/mopidy/mopidy.conf

[Install]
WantedBy=multi-user.target
