# mopidy/mopidy.service
[Unit]
Description=mopidy
After=network.target

[Service]
Type=simple
User=mopidy
ExecStart=/usr/bin/mopidy

[Install]
WantedBy=multi-user.target
