# audit/auditd.service
[Unit]
Description=Security Auditing Service

[Service]
ExecStart=/sbin/auditd -n
ExecStartPost=/sbin/auditctl -R /etc/audit/audit.rules

[Install]
WantedBy=multi-user.target
