# ircd-hybrid-serv/ircd-hybrid-serv.service
[Unit]
Description=IRC daemon services
After=network.target remote-fs.target nss-lookup.target

[Service]
Type=forking
User=ircd
PIDFile=/var/run/ircd/ircd-hybserv.pid
ExecStart=/usr/sbin/hybserv

[Install]
WantedBy=multi-user.target
