# xdm-git/xdm.service
[Unit]
Description=XDM Daemon

[Service]
PIDFile=/opt/xdm/xdm.pid
EnvironmentFile=/etc/conf.d/xdm
ExecStart=/usr/bin/python2 /opt/xdm/XDM.py --daemonize --pidfile=/opt/xdm/xdm.pid --port=${PORT} --datadir=${DIR} --config=${CONF}  --apiPort=${APIPORT}
Type=forking
User=xdm
Group=xdm

[Install]
WantedBy=multi-user.target
