# jabberd2/jabberd.service
[Unit]
Description=XMPP server
After=network.target

[Service]
User=jabber
Group=jabber
ExecStart=/usr/bin/jabberd

[Install]
WantedBy=multi-user.target
