# tlsdate/tlsdated.service
[Unit]
Description=Secure rdate

[Service]
Type=simple
GuessMainPID=yes
EnvironmentFile=/etc/conf.d/tlsdated
ExecStart=/usr/bin/tlsdated ${TLSDATED_OPTS}
ExecReload=/bin/kill -HUP $MAINPID
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
