# tlsdate/tlsdate.service
[Unit]
Description=Secure rdate

[Service]
Type=forking
GuessMainPID=yes
EnvironmentFile=/etc/conf.d/tlsdate
ExecStart=/usr/bin/tlsdated $DAEMON_OPTS
ExecReload=/bin/kill -HUP $MAINPID
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
