# open-fcoe/fcoe.service
[Unit]
Description=Open-FCoE initiator daemon
After=syslog.target network.target

[Service]
Type=simple
EnvironmentFile=/etc/conf.d/fcoe
ExecStartPre=/sbin/modprobe -qa $SUPPORTED_DRIVERS
ExecStart=/usr/sbin/fcoemon $FCOEMON_OPTS

[Install]
WantedBy=multi-user.target
Also=lldpad.socket
Also=fcoemon.socket
