# hpasm/hpasmd.service
[Unit]
Description=Proliant System Health Monitor

[Service]
Type=forking
ExecStart=/usr/sbin/hpasmd
Restart=on-abort

[Install]
WantedBy=multi-user.target
