# phpenv/php-fpm_env.service
[Unit]
Description=The PHP FastCGI Process Manager
After=syslog.target network.target

[Service]
Type=forking
PIDFile=/run/php-fpm/php-fpm_env.pid
PrivateTmp=true
TimeoutStartSec=10
ExecStart=/bin/php-fpm_env --pid /run/php-fpm/php-fpm_env.pid
ExecReload=/bin/kill -USR2 $MAINPID
[Install]
WantedBy=multi-user.target