# xe-guest-utilities/xe-linux-distribution.service
# Contributor:  Brandon Thomas <brandon.thomas@gmail.com>
[Unit]
Description=Citrix XenServer Tools
After=network.target

[Service]
Type=simple
PIDFile=/var/run/xe-daemon.pid
ExecStartPre=/usr/sbin/xe-linux-distribution /var/cache/xe-linux-distribution
ExecStartPre=/bin/mount -t xenfs none /proc/xen
ExecStart=/usr/sbin/xe-daemon -p /var/run/xe-daemon.pid
ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/bin/kill -TERM $MAINPID
ExecStop=/bin/umount /proc/xen

[Install]
WantedBy=multi-user.target
