# svxlink/svxlink.service
[Unit]
Description=Echolink Software
After=network.target
After=sound.target

[Service]
Type=forking
EnvironmentFile=/etc/default/svxlink
ExecStart=/bin/su -s /bin/bash -l -c "/usr/bin/svxlink --daemon --pidfile=${PIDFILE} --logfile=${LOGFILE} --config=${CFGFILE} --runasuser=${RUNASUSER}"
Restart=always
PIDFile=/run/svxlink.pid

[Install]
WantedBy=multi-user.target
