# tshock/tshock.service
[Unit]
Description=A Terraria Server administration modification that runs on the open source TerrariaAPI.

[Service]
ExecStart=/srv/tshock/tshock.sh start
ExecStop=/srv/tshock/tshock.sh stop
Type=forking

[Install]
WantedBy=multi-user.target

