# prelude-manager/prelude-manager.service
[Unit]
Description=High availability server that accepts secured connections from distributed sensors and/or other Managers

[Service]
Type=forking
PIDFile=/var/run/prelude-manager/prelude-manager.pid
ExecStart=/usr/bin/prelude-manager --daemon --pidfile /var/run/prelude-manager/prelude-manager.pid

[Install]
WantedBy=multi-user.target
