# rstudio-server/rstudio-server.service
[Unit]
Description=RStudio Server
After=network.target

[Service]
ExecStart=/usr/lib/rstudio-server/bin/rserver --server-daemonize=0
User=rstudio-server

[Install]
WantedBy=multi-user.target
