# postfwd/postfwd.service
[Unit]
Description=Postfwd daemon for Postfix
After=network.target
Requires=network.target

[Service]
Type=forking
User=nobody
Group=nobody
ExecStart=/usr/bin/postfwd --shortlog --summary=600 --cache=600 --cache-rbl-timeout=3600 --cleanup-requests=1200 --cleanup-rbls=1800 --cleanup-rates=1200 --daemon --file=/etc/postfix/postfwd.cf --interface=127.0.0.1 --port=10040 --pidfile=/var/tmp/postfwd.pid

[Install]
WantedBy=multi-user.target
