# dualserver/dualserverd.service
[Unit]
Description=Dual DHCP DNS Server
After=network.target

[Service]
User=root
Type=forking
ExecStart=/usr/bin/dualserverd -i /etc/dualserver/dualserver.ini -s /run/dualserver.state -l /var/log/dualserver.log
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
