# guacd/guacd.service
[Unit]
Description=Guacamole Proxy Daemon
After=syslog.target network.target

[Service]
Type=forking
PIDFile=/var/run/guacd.pid
ExecStart=/usr/sbin/guacd -p /var/run/guacd.pid

[Install]
WantedBy=multi-user.target
