# dsc-cassandra/cassandra.service
[Unit]
Description=Datastax Cassandra NoSQL database
After=syslog.target
After=network.target

[Service]
Type=forking
User=cassandra
PIDFile=/var/run/cassandra/cassandra.pid
WorkingDirectory=/var/lib/cassandra
EnvironmentFile=-/etc/conf.d/cassandra
ExecStart=/usr/bin/cassandra -p /var/run/cassandra/cassandra.pid

[Install]
WantedBy=multi-user.target
