# dionaea-git/dionaea.service
[Unit]
Description=Dionaea
After=network.target

[Service]
Type=forking
PIDFile=/var/run/dionaea.pid
ExecStart=/opt/dionaea/bin/dionaea -u dionaea -g dionaea -r /opt/dionaea -p /var/run/dionaea.pid -D
ExecStop=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target

