# autopoweroff/autopoweroff.service
[Unit]
Description=Autopoweroff Daemon

[Service]
Type=forking
ExecStart=/usr/sbin/autopoweroffd

[Install]
WantedBy=multi-user.target
