# armagetronad-dedicated/armagetronad.service
[Unit]
Description=Armagetron Advanced Dedicated Server
After=network.target

[Service]
Type=simple
User=armagetron
Group=armagetron
ExecStart=/usr/bin/armagetronad-dedicated

[Install]
WantedBy=multi-user.target
