# hpfall-git/hpfall.service
[Unit]
Description=Disk protection for HP machines.

[Service]
Type=forking
EnvironmentFile=/etc/conf.d/hpfall
ExecStart=/usr/bin/hpfall $DEVICE

[Install]
WantedBy=multi-user.target
