# statsite-git/statsite.service
[Unit]
Description=Statsite
After=network.target

[Service]
ExecStart=/usr/bin/statsite -f /etc/statsite.conf

[Install]
WantedBy=multi-user.target
