# libguestfs/guestfsd.service
[Unit]
Description=Libguestfs live service
BindTo=dev-virtio\x2dports-org.libguestfs.channel.0.device
After=dev-virtio\x2dports-org.libguestfs.channel.0.device

[Service]
ExecStart=-/usr/sbin/guestfsd -r
Restart=always
RestartSec=0

[Install]
