# ctkarch-autostartx/autostartx.install
post_install() {
  echo 'To use this:'
  echo ' configure your username in /etc/systemd/system/autostartx.service.d/user.conf'
  echo ' run `systemctl enable autostartx`'
}

post_upgrade() {
  post_install
}
