# quassel-indicator/quassel.service
[Unit]
Description=Quassel Daemon
After=network.target

[Service]
EnvironmentFile=/etc/conf.d/quassel
PIDFile=/var/run/quassel.pid
ExecStart=/bin/su - $QUASSEL_USER -c "/usr/bin/quasselcore --configdir=/var/lib/quassel"

[Install]
WantedBy=multi-user.target
