# yrmcds/yrmcds.service
[Unit]
Description=yrmcds memcached compatible KVS
ConditionPathExists=/etc/yrmcds.conf
After=network.target

[Service]
Type=simple
ExecStart=/usr/sbin/yrmcdsd
ExecReload=/bin/kill -HUP $MAINPID
ExecStartpre=/bin/mkdir -p /var/tmp/yrmcds

[Install]
WantedBy=multi-user.target

