# 7robot-eurobot-git/hubmanager.service
[Unit]
Description=Start hubmanager.

[Service]
User=nobody
Group=nobody
Type=simple
ExecStart=/usr/bin/hubmanager 1298@text 1299@text 1300@text 1301@atp 1302@atp 1303@atp 1304@atp 1305@atp 1306@atp 1307@atp 1308@atp 1309@atp

[Install]
WantedBy=multi-user.target
