# 7robot-eurobot-git/jsbot.service
[Unit]
Description=Start jsbot daemon.

[Service]
Type=simple
ExecStart=/usr/bin/jsbot

[Install]
WantedBy=multi-user.target
