# 7robot-eurobot-git/tcphub@.service
[Unit]
Description=Start a tcphub with a char decoder on specified port.

[Service]
Type=simple
ExecStart=/usr/bin/tcphub %I char

[Install]
WantedBy=multi-user.target
