# 7robot-eurobot-git/usbdaemon.service
[Unit]
Description=Start usbdaemon.
After=hubmanager.service

[Service]
Type=simple
ExecStart=/usr/bin/usbdaemon -c /etc/usbdaemon.conf -w 2000
Restart=always

[Install]
WantedBy=multi-user.target
