# 7robot-eurobot-git/usbdaemon.service
[Unit]
Description=Start usbdaemon.

[Service]
Type=simple
ExecStart=/usr/bin/usbdaemon -c /etc/usbdaemon.conf

[Install]
WantedBy=multi-user.target
