# 7robot-eurobot-git/usbdaemon@.service
[Unit]
Description=Start usbdaemon.
Conflicts=usbdaemon.service

[Service]
Type=simple
ExecStart=/usr/bin/usbdaemon -c /etc/usbdaemon.conf -h %I

[Install]
WantedBy=multi-user.target
