# i2pd/i2pd.service
[Unit]
Description=I2P Router in C++
After=network.target

[Service]
User=i2p
Type=simple
ExecStart=/usr/bin/i2p
ExecReload=/usr/bin/kill -HUP $MAINPID
KillSignal=SIGINT
LimitNOFILE=8192
PrivateDevices=yes

[Install]
WantedBy=multi-user.target
