# sonarqube/sonarqube.service
[Unit]
Description=A code quality management platform.
After=network.target

[Service]
User=sonar
Type=forking
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=sonar
ExecStart=/opt/sonarqube/bin/ARCH/sonar.sh start
ExecStop=/opt/sonarqube/bin/ARCH/sonar.sh stop
PIDFile=/opt/sonarqube/run/SonarQube.pid
OOMScoreAdjust=-1000
LimitCPU=infinity
LimitFSIZE=infinity
LimitDATA=infinity
LimitCORE=0
LimitAS=infinity
LimitLOCKS=infinity

[Install]
WantedBy=multi-user.target

