# chrome-remote-desktop/chrome-remote-desktop.service
[Unit]
Description="Chrome Remote Desktop host daemon"

[Service]
ExecStart=/opt/google/chrome-remote-desktop/chrome-remote-desktop --start
ExecStop=/opt/google/chrome-remote-desktop/chrome-remote-desktop --stop
ExecReload=/opt/google/chrome-remote-desktop/chrome-remote-desktop --reload
Restart=always
TimeoutStopSec=10

[Install]
WantedBy=default.target
