# electrum-server-git/electrum-server.service
[Unit]
Description=Electrum server
Wants=bitcoin.service
After=bitcoin.service

[Service]
Type=forking
PIDFile=/run/electrum-server.pid
User=electrum-server
Group=electrum-server
WorkingDirectory=/usr/share/electrum-server
ExecStart=/usr/bin/python2 -u server.py
ExecStop=/usr/bin/python2 server.py stop
Restart=always
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
