# electrum-server-git/electrum-server.service
[Unit]
Description=Electrum server
After=network.target

[Service]
Type=forking
PIDFile=/run/electrum-server.pid
User=electrum-server
Group=electrum-server
WorkingDirectory=/usr/share/electrum-server
ExecStart=/usr/bin/python2 -u server.py &> /var/log/electrum-server.log
ExecStop=/usr/bin/python2 server.py stop
Restart=always

[Install]
WantedBy=multi-user.target
