# fusioninventory-agent/fusioninventory-agent.service
[Unit]
Description=FusionInventory Agent
After=network.service

[Service]
ExecStart=/usr/bin/fusioninventory-agent

[Install]
WantedBy=multi-user.target
