# steam-standalone/steam-standalone.service
[Unit]
Description=starts instance of steam bigpicture using xinit
After=network.target
Conflicts=xbmc.service

[Service]
User=steam
Group=steam
Type=simple
PIDFile=%h/.steam/steam.pid
ExecStart=/usr/bin/xinit /usr/bin/steam -bigpicture -- :1 -nolisten tcp
Restart=on-failure

[Install]
WantedBy=multi-user.target
