# steam-standalone/steam-standalone.service
[Unit]
Description=standalone instance of steam bigpicture using xinit
After=network.target

[Service]
User=steam
Group=steam
Type=simple
PIDFile=%h/.steam/steam.pid
ExecStart=/usr/bin/startx
ExecStartPre=/usr/bin/sleep 1
Restart=on-failure

[Install]
WantedBy=multi-user.target
