# torque/torque-scheduler.service
[Unit]
Description=TORQUE scheduler
Wants=torque-server.service
After=torque-server.service

[Service]
Type=forking
PIDFile=/var/spool/torque/sched_priv/sched.lock
ExecStart=/usr/sbin/pbs_sched

[Install]
WantedBy=multi-user.target
