# torque/torque-node.service
[Unit]
Description=TORQUE node
Wants=basic.target
After=basic.target network.target

[Service]
Type=forking
PIDFile=/var/spool/torque/mom_priv/mom.lock
ExecStart=/usr/sbin/pbs_mom

[Install]
WantedBy=multi-user.target 
