# quarkcoin-daemon/quarkcoind@.service
[Unit]
Description=Quarkcoin Daemon
After=network.target 

[Service]
User=%i
Group=users
Type=simple
ExecStart=/usr/bin/quarkcoind

[Install]
WantedBy=multi-user.target
