# cinder/cinder-volume.service
[Unit]
Description=cinder volume server

[Service]
ExecStart=/usr/bin/cinder-volume --config-file=/etc/cinder/cinder.conf --log-file=/var/log/cinder/cinder-volume.log
User=cinder

[Install]
WantedBy=multi-user.target
