# cinder/cinder-scheduler.service
[Unit]
Description=cinder scheduler server

[Service]
ExecStart=/usr/bin/cinder-scheduler --config-file=/etc/cinder/cinder.conf --log-file=/var/log/cinder/cinder-scheduler.log
User=cinder

[Install]
WantedBy=multi-user.target
