# anope-devel/anope.service
[Unit]
Description=Anope IRC Services
Requires=network.target
After=network.target
 
[Service]
Type=simple
PIDFile=/run/anope/anope.pid
ExecStart=/usr/bin/services --confdir=/etc/anope --dbdir=/var/lib/anope --localedir=/usr/lib/anope/locales --logdir=/var/log/anope --modulesdir=/usr/lib/anope --nofork
ExecReload=/bin/kill -1 $MAINPID
User=anope
Group=anope
 
[Install]
WantedBy=multi-user.target
