# openafs/openafs-server.service
[Unit]
Description=OpenAFS Server Service
After=syslog.target network.target

[Service]
EnvironmentFile=-/etc/conf.d/openafs
ExecStart=/usr/sbin/bosserver $BOSSERVER_ARGS
ExecStop=/usr/bin/bos shutdown localhost -wait -localauth

[Install]
WantedBy=multi-user.target
