# openafs/openafs-client.service
[Unit]
Description=OpenAFS Client Service
After=syslog.target network.target

[Service]
Type=forking
EnvironmentFile=/etc/conf.d/openafs
ExecStartPre=/sbin/modprobe libafs
ExecStart=/usr/sbin/afsd $AFSD_ARGS
ExecStop=/bin/umount /afs
ExecStop=/usr/sbin/afsd -shutdown
ExecStop=/usr/bin/rmmod libafs
KillMode=none

[Install]
WantedBy=multi-user.target
