# gssproxy/gssproxy.service
[Unit]
Description=GSSAPI Proxy Daemon
# I still can't get over how ugly these are ;_;
Before=nfs-secure.service nfs-secure-server.service
Requires=proc-fs-nfsd.mount

[Service]
ExecStart=/usr/bin/gssproxy -D
Type=forking
PIDFile=/var/run/gssproxy.pid

[Install]
WantedBy=multi-user.target
