# newrelic-php/newrelic-daemon.service
[Unit]
Description=NewRelic Proxy Daemon PHP5
After=network.target

[Service]
Type=forking
PIDFile=/run/newrelic-php.pid
ExecStart=/usr/bin/newrelic-daemon -c /etc/newrelic/newrelic.cfg -p /var/run/newrelic-php.pid

[Install]
WantedBy=multi-user.target
