# cntlm-svn/cntlm-svn.service
[Unit]
Description=Cntlm-svn Proxy Server
After=network.target

[Service]
ExecStart=/usr/sbin/cntlm -f

[Install]
WantedBy=multi-user.target

