# openct-nopcsc/openct.service
[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/sh -c "mkdir -p /var/run/openct && /usr/sbin/openct-control init"
ExecStop=/usr/sbin/openct-control shutdown

[Install]
WantedBy=multi-user.target
