# upsmon/upsmon.service
[Unit]
Description=Monitor for Riello UPS devices

[Service]
Type=forking
ExecStart=/opt/upsmon/upsagent
ExecStop=killall upsagent

[Install]
WantedBy=multi-user.target
