# archiveopteryx/archiveopteryx.service
[Unit]
Description=Archiveopteryx mail server
After=network.target

[Service]
Type=forking
PIDFile=/var/lib/archiveopteryx/archiveopteryx.pid
ExecStart=/usr/bin/aox start
ExecStop=/usr/bin/aox stop

[Install]
WantedBy=multi-user.target
