# bgbilling60/bgbilling60.service
[Unit]
Description=BGBilling Server daemon
Requires=activemq.service mysqld.service
After=activemq.service mysqld.service

[Service]
EnvironmentFile=/etc/conf.d/bgbilling60
Type=forking
ExecStart=/opt/bgbilling60/server.sh start
ExecStop=/opt/bgbilling60/server.sh stop

[Install]
WantedBy=multi-user.target
