# python2-webiopi/webiopi.service
[Unit]
Description=WebIOPi daemon

[Service]
ExecStart=/usr/bin/python2 -m webiopi 8000

[Install]
WantedBy=multi-user.target
