# ivideon-video-server-nogui/videoserverd.service
[Unit]
Description=Ivideon videoserverd
After=network.target
Requires=network.target

[Service]
Type=simple
ExecStart=/opt/ivideon/ivideon-server/auto_respawn --user root --pid /var/lock/videoserverd.pid /opt/ivideon/ivideon-server/videoserver -c /root/.IvideonServer/videoserverd.config

[Install]
WantedBy=multi-user.target
