# salt-api/salt-api.service
[Unit]
Description=Salt API daemon
Requires=salt-master.service

[Service]
Type=forking
PIDFile=/run/salt-api.pid
ExecStart=/usr/bin/salt-api -d
Restart=on-failure

[Install]
WantedBy=multi-user.target
