# courier-authlib/authdaemond.service
[Unit]
Description=Courier Authentification Daemon

[Service]
Type=forking
ExecStart=/usr/sbin/authdaemond start
ExecStop=/usr/sbin/authdaemon stop
PIDFile=/run/authdaemon/pid

[Install]
WantedBy=multi-user.target
