# orientdb/orientdb.service
[Unit]
Description=OrientDB NoSQL Graph-Document DBMS

[Service]
Type=simple
User=orient
ExecStart=/opt/orientdb/bin/server.sh
ExecStop=/opt/orientdb/bin/shutdown.sh

[Install]
WantedBy=multi-user.target
