# upsource-eap/upsource.service
[Unit]
Description=Upsource server
After=network.target

[Service]
EnvironmentFile=/etc/conf.d/upsource
ExecStart=/opt/upsource-eap/bin/upsource.sh start
ExecStop=/opt/upsource-eap/bin/upsource.sh stop
Type=forking
PIDFile=/opt/upsource-eap/logs/upsource.pid
Restart=no
TimeoutSec=30
RemainAfterExit=yes
User=upsource

[Install]
WantedBy=multi-user.target

