# mediabrowser-server/mediabrowser-server.service
[Unit]
Description=Media Browser Server is a home media server built using other popular open source technologies.

[Service]
WorkingDirectory=/usr/lib/mediabrowser-server
ExecStart=/usr/bin/mediabrowser-server
User=mediabrowser
Restart=on-abort

[Unit]
After=network.target

[Install]
WantedBy=multi-user.target
