# mediabrowser-server/mediabrowser.service
[Unit]
Description=Media Server and Transcoder

[Service]
ExecStart=/usr/bin/mediabrowser-server

[Install]
WantedBy=multi-user.target
