# mediabrowser-server/mediabrowser-server.service
[Unit]
Description=Media Server and Transcoder

[Service]
WorkingDirectory=/opt/mediabrowser-server
ExecStart=/usr/bin/mediabrowser-server

[Install]
WantedBy=multi-user.target
