# havp/havp.service
[Unit]
Description=Virus Scanning Web Proxy
After=network.target

[Service]
Type=forking
PIDFile=/run/havp.pid
ExecStart=/usr/bin/havp -c /etc/havp/havp.config

[Install]
WantedBy=multi-user.target
