# v4l-mjpg-stream/v4l-mjpg-stream.service
[Unit]
Description=Stream motion JPEG from a Video4Linux enabled camera

[Service]
Type=simple
User=v4l-mjpg-stream
ExecStart=/usr/bin/mjpg_streamer -o "output_http.so -w /usr/share/mjpeg-streamer/www"

[Install]
WantedBy=basic.target
