# kiwiirc/kiwiirc.service
[Unit]
Description=Kiwi IRC web client
After=network.target

[Service]
Type=forking
PIDFile=/usr/lib/kiwiirc/kiwiirc.pid
ExecStart=/usr/bin/kiwi start
ExecReload=/usr/bin/kiwi reconfig
ExecStop=/usr/bin/kiwi stop

[Install]
WantedBy=multi-user.target

