# subvein/subvein.service
[Unit]
Description=Dedicated server for subvein
After=network.target

[Service]
ExecStart=/usr/bin/subvein-server
Type=simple
User=subvein
Group=subvein

[Install]
WantedBy=multi-user.target
