# rackspace-monitoring-agent-git/rackspace-monitoring-agent.service

[Unit]
Description=Agent for monitoring CPU and RAM of Rackspace cloud servers

[Service]
ExecStart=/usr/bin/rackspace-monitoring-agent

[Install]
WantedBy=multi-user.target
