# mfi-controller/mfi.service
[Unit]
Description=mFi Controller
After=network.target

[Service]
User=mfi
ExecStart=/usr/bin/java -jar /opt/mfi/lib/ace.jar start
ExecStop=/usr/bin/java -jar /opt/mfi/lib/ace.jar stop

[Install]
WantedBy=multi-user.target
