# udev-usbnet-exz/10-moto-ezx-usbnet.rules
#2011-03-08 Overdose: student.fmf@gmail.com
#rules for usbnet [ motorola a1200 (EZX Platform) ]

ACTION=="add", KERNEL=="usb[0-9]*", SUBSYSTEM=="net", SUBSYSTEMS=="usb", DRIVERS=="usb", ATTRS{manufacturer}=="Motorola", ATTRS{product}=="Motorola USBLAN", RUN+="/sbin/ifconfig usb%n 192.168.16.1"
ACTION=="remove", KERNEL=="usb[0-9]*", SUBSYSTEM=="net", SUBSYSTEMS=="usb", DRIVERS=="usb", ATTRS{manufacturer}=="Motorola", ATTRS{product}=="Motorola USBLAN", RUN+="/sbin/ifconfig usb%n down"
