# udpxy/udpxy.service
[Unit]
Description=UDP multicast to TCP (HTTP) relay proxy
After=network.target

[Service]
Type=forking
EnvironmentFile=/etc/conf.d/udpxy
PIDFile=/run/udpxy.pid
ExecStart=/usr/bin/udpxy $UDPXYOPTS

[Install]
WantedBy=multi-user.target
