# mydnsjpd-git/mydnsjpd.service
[Unit]
Description=MyDNS.jp DDNS record updater daemon
Requires=network.target
After=network-online.target

[Service]
Type=forking
ExecStart=/usr/bin/mydnsjpd -df /etc/mydnsjpd.conf

[Install]
WantedBy=multi-user.target
