# kwakd/kwakd.service
[Unit]
Description=A web server serving blank html pages

[Service]
Type=forking
User=kwakd
EnvironmentFile=/etc/conf.d/kwakd
ExecStart=/usr/bin/kwakd $KWAKD_ARGS
Restart=on-abort

[Install]
WantedBy=multi-user.target
