# kwakd/kwakd.service
[Unit]
Description=A web server serving blank html pages

[Service]
Type=forking
EnvironmentFile=/etc/conf.d/kwakd
ExecStart=/bin/su -s /bin/sh -c "/usr/bin/kwakd $KWAKD_ARGS" kwakd
Restart=on-abort

[Install]
WantedBy=multi-user.target
