# sickbeard/sickbeard.service
[Unit]
Description=SickBeard daemon
After=network.target

[Service]
EnvironmentFile=/etc/conf.d/sickbeard
User=sickbeard
Group=sickbeard
ExecStart=/usr/bin/sickbeard --nolaunch --daemon --quiet --config "$SB_CONF" --datadir "$SB_DATA" --port "$SB_PORT" --pidfile=/run/sickbeard/sickbeard.pid
Type=forking
TimeoutStopSec=20
PIDFile=/run/sickbeard/sickbeard.pid

[Install]
WantedBy=multi-user.target
