# dyfi-update/dyfi-update.service
[Unit]
Description=dy.fi update daemon
After=network.target

[Service]
User=dyfi-update
ExecStart=/usr/bin/dyfi-update -f /etc/dyfi-update.conf

[Install]
WantedBy=multi-user.target
