# torrent-leecher/torrent-leecher.service
[Unit]
Description=Torrent Leecher

[Service]
ExecStart=/opt/torrent-leecher/run

[Install]
WantedBy=multi-user.target
