# mediagoblin-git/mediagoblin.service
[Unit]
Description=Mediagoblin
 
[Service]
Type=forking
ExecStart=/bin/su -s /bin/sh -c '/usr/share/webapps/mediagoblin/bin/paster serve \
             /usr/share/webapps/mediagoblin/paste_local.ini \
             --pid-file=/var/run/mediagoblin/mediagoblin.pid \
             --log-file=/var/run/mediagoblin/mediagoblin.log \
             --daemon \
             --server-name=fcgi' mediagoblin
ExecStop=/bin/su -s /bin/sh -c '/usr/share/webapps/mediagoblin/bin/paster serve \
            --pid-file=/var/run/mediagoblin/mediagoblin.pid \
            /usr/share/webapps/mediagoblin/paste_local.ini stop' mediagoblin
PIDFile=/var/run/mediagoblin/mediagoblin.pid

[Install]
WantedBy=multi-user.target
