# utserver/utserver.service
[Unit]
Description=uTorrent Server Daemon
After=network.target

[Service]
Type=forking
ExecStart=/usr/bin/utserver -configfile /etc/utserver.conf -settingspath /srv/utserver/settings/ -daemon
User=utserver

[Install]
WantedBy=multi-user.target
