# randomsound/randomsound.service
[Unit]
Description=An ALSA sound card related entropy gathering daemon

[Service]
Type=forking
EnvironmentFile=/etc/conf.d/randomsound
ExecStart=/usr/sbin/randomsound -D -m ${min_entropy} -M ${max_entropy}

[Install]
WantedBy=multi-user.target
