# lxc-docker-git/docker.service
[Unit]
Description=Docker Daemon
After=network.target

[Service]
Type=simple
ExecStartPre=/bin/mount --make-rprivate /
ExecStart=/usr/sbin/docker -d

[Install]
WantedBy=multi-user.target
