# graylog2-server/graylog2.service
[Unit]
Description=Graylog2 management server
After=network.target mongodb.service
Requires=mongodb.service

[Service]
User=nobody
PIDFile=/run/graylog2-server/graylog2.pid
ExecStart=/usr/bin/java -jar /usr/lib/graylog2-server.jar -f /etc/graylog2.conf -p /run/graylog2-server/graylog2.pid

[Install]
WantedBy=multi-user.target
