# graylog2-server/graylog2.service
[Unit]
Description=Graylog2 log management server
After=syslog.target network.target

[Service]
Type=forking
Environment=GRAYLOG2_CONF=/etc/graylog2.conf
ExecStart=/usr/bin/java -jar /usr/lib/graylog2-server.jar -f ${GRAYLOG2_CONF}

[Install]
WantedBy=multi-user.target