# graylog2-server/graylog2.service
[Unit]
Description=Graylog2 log management server
After=syslog.target network.target

[Service]
Environment=JAVA_HOME=/usr/lib/jvm/java-7-openjdk
Environment=GRAYLOG2_CONF=/etc/graylog2.conf
PIDFile=/run/graylog2.pid
ExecStart=/usr/bin/graylog2-server -f ${GRAYLOG2_CONF} -p /run/graylog2.pid

[Install]
WantedBy=multi-user.target