# headphones/headphones.service
[Unit]
Description=Headphones Daemon

[Service]
EnvironmentFile=/etc/conf.d/headphones
ExecStart=/usr/bin/env python2 /opt/headphones/Headphones.py --config ${HP_CONF} --datadir ${HP_DIR} --port ${HP_PORT}
Type=simple
User=headphones
Group=headphones

[Install]
WantedBy=multi-user.target
