# headphones/headphones.service
[Unit]
Description=Headphones Daemon

[Service]
EnvironmentFile=/etc/conf.d/headphones
ExecStart=/opt/headphones/Headphones.py --config ${HP_CONF} --datadir ${HP_DIR} --port ${HP_PORT} --nolaunch
Type=simple
User=headphones
Group=downloads

[Install]
WantedBy=multi-user.target
