# smack-system/smack.mount
[Unit]
Description=smack filesystem
Before=local-fs.target

[Mount]
What=smackfs
Where=/smack
Type=smackfs
Options=smackfsdef=*
