# smack-system/smack.service
[Unit]
Description=Load SMACK rules
Requires=smack.mount
After=smack.mount

[Service]
Type=oneshot
# use this script for now
ExecStart=/etc/rc.d/smack start
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
