# shrew-vpn-client-unstable/iked.service
[Unit]
Description=Internet Key Exchange daemon

[Service]
Type=forking
ExecStart=/usr/sbin/iked

[Install]
WantedBy=multi-user.target
