# p2pool/p2pool.service
[Unit]
Description=P2Pool Daemon
After=network.target

[Service]
User=p2pool
Group=p2pool
EnvironmentFile=/etc/conf.d/p2pool
ExecStart=/usr/bin/p2poold --datadir=/var/lib/p2pool $ARGS

[Install]
WantedBy=multi-user.target
