# keysid/keysid.service
[Unit]
Description=Linux keys intercepting daemon

[Service]
Type=forking
ExecStart=/usr/sbin/keysid

[Install]
WantedBy=multi-user.target
