# elsa/elsa.cron
* * * * * root /usr/bin/perl /usr/share/elsa/web/cron.pl -c /etc/elsa_web.conf > /dev/null 2>&1
