# chromiumos-power-manager/useractd.service
[Unit]
Description=User Activity Monitor for ChromeOS Power Daemon
Requires=dbus.service powerd.service
After=dbus.service powerd.service

[Service]
Type=simple
ExecStart=/bin/sh -c '/usr/bin/useractd > /var/log/useractd.log 2>&1'

[Install]
WantedBy=multi-user.target
