# ffdecsawrapper-git/ffdecsawrapper.service
[Unit]
Description=FFdecsawrapper - FFdecsa Empowered Softcam
After=oscam.service 
Requires=oscam.service

[Service]
Type=forking
Nice=-3
User=root
PIDFile=/run/ffdecsawrapper.pid
EnvironmentFile=/etc/conf.d/ffdecsawrapper
ExecStartPre=/usr/bin/modprobe dvbloopback $DVBLOOPBACK_ARGS
TimeoutSec=2
ExecStart=/usr/bin/ffdecsawrapper $FFDECSAWRAPPER_ARGS --daemon --pidfile /run/ffdecsawrapper.pid --cam-dir $CAMDIR --log $LOGFILE

[Install]
WantedBy=multi-user.target
