# cjdns-git/cjdns.service
[Unit]
Description=A routing engine designed for security, scalability, speed and ease of use

[Service]
Type=forking
GuessMainPID=yes
EnvironmentFile=/etc/conf.d/cjdns
ExecStart=/usr/lib/systemd/scripts/cjdns $CJDNS_CONFIG
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
