# klnagent/klnagent.service
[Unit]
Description=Kaspersky Network Agent
After=network.target
After=systemd-journald.service

[Service]
Type=forking
PIDFile=/var/run/klnagent.pid
ExecStart=/opt/kaspersky/klnagent/sbin/klnagent
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=always
RestartSec=300

[Install]
WantedBy=multi-user.target
