# g9led/g9ledd.service
[Unit]
Description=Logitech G9 mouse LED color-changing daemon by time interval

[Service]
Type=forking
ExecStart=/bin/sh -c '/usr/bin/g9ledd'
 
[Install]
WantedBy=multi-user.target