# openvpn_proxy/openvpn_proxy.service
[Unit]
Description=Convert OpenVPN TCP connections into UDP connections
After=network.target

[Service]
Type=forking
EnvironmentFile=-/etc/conf.d/openvpn_proxy
ExecStart=/usr/bin/openvpn_proxy -p $LISTEN_PORT -H $LISTEN_ADDRESS -s $OPENVPN_PORT -S $OPENVPN_ADDRESS

[Install]
WantedBy=multi-user.target
