# inadyn/inadyn.service
[Unit]
Description=Dynamic DNS updater
After=network.target
Documentation=man:inadyn(8)

[Service]
Type=forking
ExecStart=/usr/bin/inadyn

[Install]
WantedBy=multi-user.target
