# buddycloud-server-git/buddycloud-server.service
[Unit]
Description=buddycloud server (git version)
After=network.target prosody.service ejabberd.service openfire.service

[Service]
User=buddycloud
WorkingDirectory=/opt/buddycloud-server-git
Type=simple
ExecStart=/usr/bin/node -- lib/main.js

[Install]
WantedBy=multi-user.target
