# 42coin-daemon/42coind@.service
[Unit]
Description=42 Coin Daemon
After=network.target 

[Service]
User=%i
Group=users
Type=simple
ExecStart=/usr/bin/42coind

[Install]
WantedBy=multi-user.target
