# iodined-systemd/iodined.service
[Unit]
After=network.target
Requires=network.target

[Service]
EnvironmentFile=/etc/conf.d/iodined
ExecStart=/usr/sbin/iodined -f $TUN_IP $TOP_DOMAIN

[Install]
WantedBy=multi-user.target
