# mldonkey-cvs/mldonkey.service
[Unit]
Description=MLDonkey - A multi-network P2P client
After=network.target

[Service]
Type=forking
#User=
ExecStart=/bin/sh -c "mlnet &"
#Restart=always
#RestartSec=1

[Install]
WantedBy=multi-user.target
