# xen/xend.service
[Unit]
Description=Xen common services

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/etc/rc.d/xencommons start
ExecStop=/etc/rc.d/xencommons stop

[Install]
WantedBy=multi-user.target 
