# heat/heat-api-cloudwatch.service
[Unit]
Description=Heat API Cloudwatch server

[Service]
ExecStart=/usr/bin/heat-api-cloudwatch
User=heat

[Install]
WantedBy=multi-user.target
