# heat/heat-api-cfn.service
[Unit]
Description=Heat API CFN server

[Service]
ExecStart=/usr/bin/heat-api-cfn
User=heat

[Install]
WantedBy=multi-user.target
