# heat/heat-api.service
[Unit]
Description=Heat API server

[Service]
ExecStart=/usr/bin/heat-api
User=heat

[Install]
WantedBy=multi-user.target
