# heat/heat-engine.service
[Unit]
Description=Heat engine server

[Service]
ExecStart=/usr/bin/heat-engine
User=heat

[Install]
WantedBy=multi-user.target
