# inputlirc/inputlircd.service
[Unit]
Description=Zeroconf LIRC daemon that reads from /dev/input/event devices
After=network.target lirc.service

[Service]
Type=forking
EnvironmentFile=/etc/inputlircd.conf
ExecStartPre=/usr/bin/mkdir -p /var/run/lirc
ExecStart=/usr/sbin/inputlircd $OPTIONS $EVENT_DEVICE
ExecStatPost=[ -S /var/run/lirc/lircd ] && ln -sf /var/run/lirc/lircd /dev/lircd

[Install]
WantedBy=multi-user.target

