# nzbget-systemd/nzbget.service
[Unit]
Description=NZBGet Daemon
DaemonAfter=network.target
 
[Service]
User=nzbget
Group=nzbget
ExecStart=/usr/bin/nzbget -c /var/lib/nzbget/nzbget.conf -D
ExecStop=/usr/bin/nzbget -Q
ExecReload=/usr/bin/nzbget -O
 
[Install]
WantedBy=multi-user.target
