# clevo-airplane-mode/60-clevo-airplane-mode.rules
KERNEL=="rfkill[0-9]*", ACTION=="change", RUN+="/usr/bin/bash /usr/bin/clevo-airplane-mode-led-control $number"
