# wolframe/wolframed.service
[Unit]
Description=Wolframe Daemon
Documentation=man:wolframed(8) man:wolframe.conf(5) http://www.wolframe.net
After=network.target
After=syslog.target

[Service]
ExecStart=/usr/sbin/wolframed -c /etc/wolframe/wolframe.conf
Type=forking
PIDFile=/var/run/wolframe/wolframed.pid

[Install]
WantedBy=multi-user.target
