# chiliproject-standalone/chiliproject.service
[Unit]
Description=Chiliproject webserver
After=network.target

[Service]
User=chiliproject
ExecStart=/usr/bin/ruby /var/lib/chiliproject/script/server -e production &>>/var/lib/chiliproject/log/production.log
ExecStop=/bin/kill -INT $MAINPID

[Install]
WantedBy=multi-user.target
