# p910nd/p910nd.service
[Unit]
Description=A small printer daemon intended for diskless workstations that passes jobs directly to the printer
After=network.target

[Service]
Type=forking
ExecStart=/usr/bin/p910nd

[Install]
WantedBy=multi-user.target
