# mkinitcpio-sd-resume/systemd-resume-setup.service
[Unit]
Description=Setup Resume Device
DefaultDependencies=no
Conflicts=shutdown.target
After=sysinit.target
Before=basic.target shutdown.target
;ConditionKernelCommandLine="resume" Currently doesn't work
ConditionPathIsReadWrite=/sys
;ConditionDirectoryNotEmpty=/dev/disk Also doesn't work

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/lib/systemd/scripts/resume-setup

