# vmware-view-client/vmware-view-usbd.service
[Unit]
Description=VMware View USB daemon

[Service]
Type=forking
ExecStart=/usr/lib/systemd/scripts/vmware-view-usbd

[Install]
WantedBy=multi-user.target
