# palm-novacom-git/novacomd.service
[Unit]
Description=Novacom WebOS daemon

[Service]
EnvironmentFile=/etc/conf.d/novacomd
ExecStart=/usr/bin/novacomd  $NOVACOMD_ARGS
