# openmoko-usb-networking/z60_openmoko_net.rules
ACTION=="add", SUBSYSTEM=="net", ATTRS{idVendor}=="1457", ATTRS{idProduct}=="5122", RUN+="/etc/rc.d/openmoko-usb-networking"
ACTION=="remove", SUBSYSTEM=="net", INTERFACE="eth1", RUN+="/etc/rc.d/openmoko-usb-networking"

