# shishi/shishid.service
[Unit]
Description=Shishi KDC Server

[Service]
ExecStart=/usr/sbin/shishid --quiet
SyslogFacility=auth

[Install]
WantedBy=multi-user.target
