# btsync-autoconfig/btsync-autoconfig@.service
[Unit]
Description=On-demand BitTorrent Sync config creator for %i
Before=btsync@.service

[Service]
User=%i
Type=oneshot
ExecStart=/bin/bash /usr/share/btsync-autoconfig/btsync-autoconfig.sh

[Install]
WantedBy=btsync@.service
