# hiawatha/php-fcgi.service
[Unit]
Description=PHP FCGI Daemon for the Hiawatha Web Server
After=network.target

[Service]
Type=forking
SyslogIdentifier=php-fcgi
ExecStart= /usr/sbin/php-fcgi
ExecStop=/usr/sbin/php-fcgi -k

[Install]
WantedBy=multi-user.target
