# counterpartyd-git/counterpartyd.service
[Unit]
Description=Counterpartyd server

[Service]
User=counterparty
Group=counterparty
ExecStart=/usr/bin/counterpartyd \
                   --config-file='/etc/counterparty/counterpartyd.conf'        \
                   --data-dir='/usr/share/counterparty'                        \
                   --log-file='/var/log/counterparty/counterpartyd-server.log' \
                   --verbose --force server

[Install]
WantedBy=multi-user.target
