# nzbmegasearch-git/nzbmegasearch.service
[Unit]
Description=NZBmegasearcH
After=network.target

[Service]
User=nzbmegasearch
Group=nzbmegasearch
ExecStart=/usr/bin/mega2
Type=simple
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
