# systemd-granola/granola.service
[Unit]
Description=Granola Core Daemon
After=systemd-user-sessions.service

[Service]
ExecStart=/usr/sbin/granola

[Install]
WantedBy=granola-connect.service
