# vmware-usbarbitrator/vmware-usbarbitrator.service
[Unit]
Description=VMware USB Arbitrator
[Service]
ExecStart=/etc/init.d/vmware-USBArbitrator start
ExecStop=/etc/init.d/vmware-USBArbitrator stop
PIDFile=/var/lock/subsys/vmware
TimeoutSec=0
RemainAfterExit=yes
[Install]
WantedBy=multi-user.target
