# riverrun-git/riverrun.service
[Unit]
Description=A minimalist book library web interface.

[Service]
ExecStart=/usr/bin/riverrun -c /etc/riverrun.conf

[Install]
WantedBy=multi-user.target

