# wolframe-git/wolframed.service
[Unit]
Description=Wolframe Daemon
After=network.target
After=syslog.target

[Service]
ExecStart=/usr/sbin/wolframed -c /etc/wolframe/wolframe.conf
Type=forking
PIDFile=/var/run/wolframe/wolframed.pid

[Install]
WantedBy=multi-user.target
