# zemberek-server/zemberek-server.service
[Unit]
Description=A Turkish spell checker server based on Zemberek NLP library
Requires=dbus.service
After=dbus.service

[Service]
Type=forking
EnvironmentFile=-/etc/conf.d/zemberek
ExecStart=/usr/bin/zemberek-server start
ExecStop=/usr/bin/zemberek-server stop
ExecReload=/usr/bin/zemberek-server restart

[Install]
WantedBy=multi-user.target