# pianod/pianod.service
[Unit]
Description=Pandora Client Daemon
After=syslog.target

[Service]
ExecStart=/usr/sbin/pianod -g audio
Restart=on-abort

[Install]
WantedBy=multi-user.target
