# tp-hdd-led/tp-hdd-led.service
[Unit]
Description=Enable ThinkPad LED as HDD LED
DefaultDependencies=true
Wants=graphical.target

[Service]
ExecStart=/usr/share/tp-hdd-led/tp-hdd-led.py

[Install]
WantedBy=multi-user.target

