# p-rout-git/p-rout-dump.service
[Unit]
Description=Dump p-rout SQL data

[Service]
Type=oneshot
ExecStart=/home/bertb/p-rout/p-rout-dump.sh
Nice=19
IOSchedulingClass=best-effort
IOSchedulingPriority=7
