# sac_linux/SACSrv.service
[Unit]
Description=eToken service aladdin-rd.ru
Documentation=/usr/share/doc/eToken

[Service]
ExecStart=/usr/bin/SACSrv
ExecReload=/bin/kill -HUP $MAINPID
StandardOutput=null

[Install]
WantedBy=multi-user.target
Alias=etoken.service
