# opennebula-unstable/opennebula.service
[Unit]
Description=OpenNebula management

[Service]
Type=oneshot
ExecStart=/bin/su - oneadmin /usr/bin/one start
ExecStop=/bin/su - oneadmin /usr/bin/one stop
TimeoutSec=0
StandardInput=tty
RemainAfterExit=yes
