# pgl-cli/pgl-update.service
[Unit]
Description=Update PeerGuardian Linux blocklists
Wants=network.target

[Service]
ExecStart=/usr/bin/pglcmd update

[Install]
WantedBy=multi-user.target
