# pgl-cli/pgl-update.timer
[Unit]
Description=Update PeerGuardian Linux blocklists every day

[Install]
WantedBy=multi-user.target

[Timer]
OnUnitActiveSec=1d
