# pgl-cli/pgl.service
[Unit]
Description=PeerGuardian Linux - an IP Blocker
Wants=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/pglcmd start
ExecStop=/usr/bin/pglcmd stop_quick
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
