# wide-dhcpv6/dhcp6c@.service
[Unit]
Description=WIDE-DHCPv6 Client on interface %I
After=network.target

[Service]
ExecStart=/usr/bin/dhcp6c -c /etc/wide-dhcpv6/dhcp6c.conf -fp /var/run/dhcp6c.pid %I

[Install]
WantedBy=multi-user.target
