# fluxfonts/fluxfonts.service
[Unit]
Description=Continual random font generator for increased privacy

[Service]
Type=forking
ExecStart=/usr/bin/fluxfontd

[Install]
WantedBy=multi-user.target

