# suricata-pfring-git/suricata.service
[Unit]
Description=Suricata IDS/IPS
After=network.target

[Service]
Type=forking
PIDFile=/run/suricata/suricata.pid

ExecStart=/sbin/suricata $OPTIONS -D
EnvironmentFile=-/etc/default/suricata

ExecReload=/usr/bin/kill -USR2 $MAINPID
KillSignal=SIGQUIT
KillMode=mixed

[Install]
WantedBy=multi-user.target

