# monitorix/monitorix.service
[Unit]
Description=A lightweight system monitoring tool that uses rrd databases.
Wants=local-fs.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/etc/rc.d/monitorix start
ExecStop=/etc/rc.d/monitorix stop

[Install]
WantedBy=multi-user.target
