# monitorix/monitorix.service
[Unit]
Description=A lightweight system monitoring tool that uses rrd databases.
Wants=local-fs.target

[Service]
Type=forking
ExecStart=/usr/bin/monitorix -c /etc/monitorix.conf -p /run/monitorix.pid
PIDFile=/run/monitorix.pid

[Install]
WantedBy=multi-user.target
