# opscenter/opscenter.service
[Unit]
Description=Visual management and monitoring solution for Cassandra
After=syslog.target

[Service]
Type=forking
PIDFile=/var/run/opscenter.pid
EnvironmentFile=-/etc/conf.d/opscenter
ExecStart=/usr/bin/opscenter -p /var/run/opscenter.pid

[Install]
WantedBy=multi-user.target
