# fgms-git-stable/fgms.service
#FGMS systemd unit contributed by ImNtReal

[Unit]
Description=FlightGear Multiplayer Server

[Service]
ExecStart=/usr/sbin/fgms

[Install]
WantedBy=multi-user.target
