# volctl2/volctld.service
[Unit]
Description=A volume control daemon

[Service]
ExecStart=/usr/bin/volctld

[Install]
WantedBy=desktop.target
