# redirex/redirex.service
[Unit]
Description=Redirect HTTP Requests to New Web Server
After=syslog.target network.target

[Service]
Type=simple
ExecStart=/usr/sbin/redirex -c /etc/redirex.conf

[Install]
WantedBy=multi-user.target
