# discretemon/89-discretemon.rules
# Hotplug detection of external monitor to be discrete (/usr/bin/discretemon)
# Resource: http://tinyurl.com/afl9mm8

ACTION=="change", SUBSYSTEM=="drm", HOTPLUG=="1", RUN+="/usr/bin/discretemon"
