# discretemon/89-discretemon.rules
# Hotplug detection of external monitor to be discrete (/usr/bin/discretemon)
# Resource: http://tinyurl.com/afl9mm8

KERNEL=="card0", SUBSYSTEM=="drm", ACTION=="change", RUN+="/usr/bin/xuserrun /usr/bin/discretemon"
