# discretemon/discretemon-resume.service
[Unit]
Description=set external monitor as discrete following resume
After=suspend.target

[Service]
User=root
Type=forking
ExecStart=/usr/bin/xuserrun /usr/bin/discretemon

[Install]
WantedBy=suspend.target
