# foldingathome-v7/foldingathome-v7.service
[Unit]
Description=Folding@home V7 distributed computing client
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/fah-v7
ExecStart=/opt/fah-v7/FAHClient --config /opt/fah-v7/config.xml --exec-directory=/opt/fah-v7 --data-directory=/opt/fah-v7

[Install]
WantedBy=multi-user.target

