# python2-fileconveyor-git/fileconveyor.service
[Unit]
Description=Transport files to a server using inotify

[Service]
ExecStart=/usr/bin/python2 -m fileconveyor.arbitrator

[Install]
WantedBy=multi-user.target
