# cpupowerd/cpupowerd.service
[Unit]
Description=cpupowerd service - undervolting vcore with Linux

[Service]
Type=forking
ExecStart=/usr/sbin/cpupowerd -d -c /etc/cpupowerd.conf

[Install]
WantedBy=multi-user.target
