# influxdb-collectd-proxy-git/influxdb-collectd-proxy.service
[Unit]
Description=Proxy that forwards collectd data to influxdb

[Service]
Type=simple
EnvironmentFile=/etc/conf.d/influxdb-collectd-proxy
ExecStart=/usr/bin/influxdb-collectd-proxy --database=$database --username=$username --password=$password --typesdb=$typesdb
User=collectd-proxy
Group=collectd-proxy

[Install]
RequiredBy=collectd.service
