# phxd-hg/phxd.service
[Unit]
Description=Python Hotline X Daemon
After=network.target

[Service]
ExecStart=/usr/bin/python2 /srv/phxd/phxd.py ; PHXD_PID="pgrep -n python"
Type=oneshot
RemainAfterExit=yes
ExecStop=pkill $PHXD_PID

[Install]
WantedBy=multi-user.target
