# lirc-zotac-git/lircd.service
[Unit]
Description=Linux Infrared Remote Control Service
ConditionPathExists=/dev/remote

[Service]
Type=forking
ExecStart=/usr/sbin/lircd --pidfile /run/lircd.pid -d /dev/remote -H zotac
PIDFILE=/run/lircd.pid

[Install]
WantedBy=multi-user.target

