# openerp/openerp.service
[Unit]
Description=Advanced OpenSource ERP and CRM server
Requires=postgresql.service
After=postgresql.service

[Service]
EnvironmentFile=/etc/conf.d/openerp-server
User=openerp
SyslogIdentifier=openerp
PIDFile=/run/openerp/openerp-server.pid
ExecStart= /usr/bin/openerp-server -c /etc/openerp/openerp-server.conf --syslog --pidfile=/run/openerp/openerp-server.pid $OPENERP_ARGS
ExecStop= /bin/kill -15 $MAINPID

[Install]
WantedBy=multi-user.target
