# soundmodem/soundmodem.service
[Unit]
Description=Soundmodem TNC

[Service]
ExecStart=/usr/sbin/soundmodem

[Install]
WantedBy=multi-user.target
