# chinadns-c/chinadns.service
[Unit]
Description=ChinaDNS Service
After=network.target

[Service]
Type=simple
User=nobody
ExecStart=/usr/bin/chinadns -l /etc/iplist.txt -c /etc/chnroute.txt -p 5353

[Install]
WantedBy=multi-user.target

