# bitrated-git/bitrated.service
[Unit]
Description=Bitrated Web service
Requires=mongodb.service
After=mongodb.service
After=network.target

[Service]
Type=simple
PIDFile=/run/bitrated.pid
User=bitrated
Group=bitrated
ExecStart=/bin/bash /opt/bitrated/bin/start.sh &
ExecStop=/usr/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
