# jetty8/jetty.service
[Unit]
Description=Jetty Web Application Server
After=syslog.target network.target

[Service]
Type=forking
PIDFile=/var/run/jetty.pid
ExecStart=/usr/bin/jetty start
ExecStop=/usr/bin/jetty stop

[Install]
WantedBy=multi-user.target

