# pattern-device-git/patterndev.service
[Unit]
Description=Pattern Device
Requires=udev-settle.service system-modules-load.service
After=udev-settle.service system-modules-load.service

[Service]
Type=oneshot
ExecStart=/etc/rc.d/patterndev start
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
