# openonload/openonload.service
[Unit]
Description=OpenOnload Service
After=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
PrivateTmp=true
ExecStart=/etc/init.d/openonload start
ExecStop=/etc/init.d/openonload stop

[Install]
WantedBy=multi-user.target
