# cloudprint-server/cloudprint-server.service
[Unit]
Description=Google Cloud Print Printer Server
After=network.target

[Service]
User=root
ExecStart=/usr/bin/chromium --type=service --enable-cloud-print-proxy --no-service-autorun --noerrdialogs --user-data-dir=/etc/cloudprint-server/
Restart=always

[Install]
WantedBy=multi-user.target
