# ss5/ss5.service
[Unit]
Description=SS5 Socks Server

[Service]
Type=forking
PIDFile=/var/run/ss5/ss5.pid
EnvironmentFile=/etc/conf.d/ss5
ExecStartPre=/usr/bin/mkdir -p /var/run/ss5
ExecStart=/bin/sh -c "/usr/sbin/ss5 -s -t"
ExecStop=/usr/bin/pkill ss5
ExecStopPost=/usr/bin/rm -f /var/run/ss5/ss5.pid

[Install]
WantedBy=multi-user.target
