# voipmonitor/voipmonitor.service
[Unit]
Description=Monitor voice quality of SIP/RTP/RTCP
Requires=network.target remote-fs.target mysqld.service
After=network.target remote-fs.target mysqld.service

[Service]
Type=forking
PIDFile=/run/voipmonitorinitscript.pid
ExecStart=/usr/bin/voipmonitor --config-file /etc/voipmonitor.conf --pid-file /run/voipmonitorinitscript.pid -v 1
SuccessExitStatus=0 1 6
