# maradns2/zoneserver.service
[Unit]
Description=Zone AXFR service for MaraDNS server
Documentation=/usr/share/doc/maradns

[Service]
ExecStart=/usr/sbin/zoneserver -f /etc/mararc
ExecReload=/bin/kill -HUP $MAINPID
StandardOutput=null

[Install]
WantedBy=multi-user.target
Alias=Zoneserver.service
