# sonar/sonar.service
[Unit]
Description=Sonar 3.x
After=network.target

[Service]
ExecStart=/opt/sonar/bin/sonar.sh start
ExecStop=/opt/sonar/bin/sonar.sh stop
ExecReload=/opt/sonar/bin/sonar.sh restart
PIDFile=/opt/sonar/run/sonar.pid
Type=forking
User=sonar

[Install]
WantedBy=multi-user.target

