# syncthing-discosrv/syncthing-discosrv.service
[Unit]
Description=Syncthing discovery server
After=network.target

[Service]
User=syncthing
Group=syncthing
ExecStart=/usr/bin/syncthing-discosrv
Restart=on-failure

[Install]
WantedBy=multi-user.target
