# python2-sentry/sentry.service
[Unit]
Description=Powerful event logger to assist in debugging and development.
After=network.target

[Service]
User=sentry
Restart=always
ExecStart=/usr/bin/sentry --config=/etc/sentry/sentry.conf.py start

[Install]
WantedBy=multi-user.target
