# ec2arch-systemd/ec2-ssh-keys.service
[Unit]
Description=Install EC2 ssh keys for root login

[Service]
ExecStart=/usr/sbin/ec2-ssh-keys

[Install]
WantedBy=multi-user.target
