# craftbukkit/craftbukkit.service
[Unit]
Description=Minecraft server implementing the Bukkit API

[Service]
ExecStart=/srv/craftbukkit/craftbukkit.sh
Type=forking
User=craftbukkit

[Install]
WantedBy=multi-user.target
