# minimserver/minimserver@.service
[Unit]
Description=MinimServer

[Service]
User=%I
ExecStart=/bin/bash /opt/minimserver/bin/startc
ExecStop=/bin/bash /opt/minimserver/bin/stopall

[Install]
WantedBy=network.target