# riak-systemd-unit/riak.service
[Unit]
Description=Riak Distributed Database

[Service]
Type=simple
User=riak
SyslogIdentifier=riak

ExecStart=  /opt/riak/bin/riak start
ExecStop=   /opt/riak/bin/riak stop

[Install]
WantedBy=multi-user.target
