# scirius-git/scirius.service
[Unit]
Description=Scirius Rule Managemant for Suricata
After= network.target
Before= nginx.service

[Service]
User=www
ExecStart=/usr/bin/python2 /usr/share/webapps/scirius/manage.py runserver
StandardOutput=null

[Install]
WantedBy=multi-user.target
