# scirius-git/scirius-fcgi.service
[Unit]
Description=Scirius Rule Managemant for Suricata
After= network.target
Before= nginx.service

[Service]
PIDFile=/var/run/scirius/scirius.pid
Type=forking
ExecStart=/usr/bin/python2 /usr/share/webapps/scirius/manage.py runfcgi socket=/var/run/scirius/fcgi.sock methode=prefork daemonize=true pidfile=/var/run/scirius/scirius.pid


[Install]
WantedBy=multi-user.target
