# scirius-git/suri-reloader.service
[Unit]
Description=Reloads Suricata when scirius.reload file is created
After=suricata.service

[Service]
Type=idle
ExecStart=/usr/bin/suri_reloader $OPTIONS
EnvironmentFile=-/etc/default/suri-reloader


[Install]
WantedBy=multi-user.target

