# systemd-journald-volatile/systemd-journald-volatile.service
[Unit]
Description=Remove Journal path for journald
DefaultDependencies=no
Before=shutdown.target reboot.target halt.target

[Service]
ExecStart=/sbin/systemd-journald-volatile-rm.sh

[Install]
WantedBy=shutdown.target reboot.target halt.target
