# aksusbd/aksusbd.service

[Unit]
Description=AKSUSB Daemon

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/sbin/aksusbd ; \
          /usr/sbin/winehasp ; \
          /usr/sbin/hasplmd -s
ExecStop=/usr/bin/killall aksusbd ; \
         /usr/bin/killall winehasp ; \
         /usr/bin/killall hasplmd
Restart=on-abort

[Install]
WantedBy=multi-user.target

