# paste-binouse/paste-binoused.service
[Unit]
Description=Paste Binouse sever

[Service]
Type=simple
SyslogIdentifier=paste-binoused
ExecStart=/usr/bin/paste-binoused -uid paste-binouse -gid paste-binouse -log-level info -capacity 4194304 -max-clients 8 -max-paste-size 131072 -ms-sleep 100 -root /srv/paste-binouse -port 11013
