# smartladder-git/smartladder.service
[Unit]
Description=SmartLadder Service
After=network.target

[Service]
Type=simple
User=nobody
ExecStart=/usr/share/smartladder/goagent+/proxy.py

[Install]
WantedBy=multi-user.target

