# cloud9-git/cloud9.service
[Unit]
Description=Cloud9 Ide

[Service]
Type=forking
User=cloud9
WorkingDirectory=/var/lib/cloud9/
ExecStart=/usr/bin/cloud9
#ExecStop=/usr/bin/openbazaar-stop

[Install]
WantedBy=multi-user.target
