# lenovo-t420-keymap/81-lenovo-keymap.rules
SUBSYSTEM=="input", ATTRS{name}=="AT Translated Set 2 keyboard", RUN+="/lib/udev/keymap input/%k lenovo-thinkpad_t420"
