# neutron/neutron-plugin-openflow-agent.service
[Unit]
Description=Neutron Openflow Plugin agent

[Service]
ExecStart=/usr/bin/neutron-ofagent-agent --config-file=/etc/neutron/neutron.conf --config-file=/etc/neutron/ml2/ml2_conf_ofa.ini --log-file=/var/log/neutron/openflow-agent.log
User=neutron

[Install]
WantedBy=multi-user.target
