# neutron/neutron-ovs-cleanup.service
[Unit]
Description=Neutron OVS cleanup

[Service]
ExecStart=/usr/bin/neutron-ovs-cleanup --config-file=/etc/neutron/neutron.conf --log-file=/var/log/neutron/ovs-cleanup.log --verbose
User=neutron

[Install]
WantedBy=multi-user.target
