# opensm/opensm-genconfig@.service
[Unit]
Description=InfiniBand OpenSM configuration generation
ConditionPathExists=|!/etc/opensm/%I.conf

[Service]
ExecStart=/usr/sbin/opensm --guid %I --create-config "/etc/opensm/%I.conf" --logfile "/var/log/opensm-%I.log"
Type=oneshot
RemainAfterExit=yes
