# opensm/opensm@.service
[Unit]
Description=Configured InfiniBand OpenSM
Before=network.target
After=opensm-genconfig@%I.service

[Service]
ExecStart=/usr/sbin/opensm --daemon --config "/etc/opensm/%I.conf"
Type=forking
