# opensm/opensm.service
[Unit]
Description=InfiniBand OpenSM
Before=network.target

[Service]
ExecStart=/usr/sbin/opensm --daemon
Type=forking
