# mcproxy-git/mcproxy.service
[Unit]
Description=mcproxy multicast proxy for IGMP/MLD (RFC 4605)
Documentation=http://mcproxy.realmv6.org/wiki/Documentation
After=network.target

[Service]
ExecStart=/usr/bin/mcproxy -f /etc/mcproxy.conf

[Install]
WantedBy=multi-user.target
