# wallproxy/wallproxy.service
[Unit]
Description=WallProxy Service
After=network.target

[Service]
Type=simple
User=nobody
ExecStart=/usr/share/wallproxy/local/startup.py

[Install]
WantedBy=multi-user.target
