# viacrypt-git/viacrypt.service
[Unit]
Description=Viacrypt server
Requires=mongodb.service redis.service
After=mongodb.service redis.service

[Service]
User=viacrypt
Group=viacrypt
WorkingDirectory=/usr/share/viacrypt
ExecStart=/usr/bin/node server.js
ExecReload=/usr/bin/kill -HUP $MAINPID
Restart=always
RestartSec=60
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
