# opendmarc/opendmarc.service
[Unit]
Description=OpenDMARC

[Service]
Type=simple
ExecStart=/usr/bin/opendmarc -f -c /etc/opendmarc/opendmarc.conf
Restart=always

[Install]
WantedBy=multi-user.target
