# libct-git/libctd.socket
[Unit]
Description = libct socket for containers management
Requires = libctd.socket

[Socket]
ListenStream = /run/libctd.socket

[Install]
WantedBy = sockets.target
Also = libctd.socket
