# libct-git/libctd.service
[Unit]
Description = libct daemon for containers management

[Service]
Type = simple
PIDFile = /run/libctd.pid
ExecStart = /usr/sbin/libctd -s /run/libctd.socket -P /run/libctd.pid

[Install]
WantedBy = multi-user.target
