# firehol/firehol.service
[Unit]
Description=Firehol stateful packet filtering firewall builder

[Service]
Type=oneshot
ExecStart=/usr/sbin/firehol start
ExecStop=/usr/sbin/firehol stop
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
