# systemd-hiawatha-units/hiawatha.service
[Unit]
Description=Secure and advanced webserver
After=network.target

[Service]
Type=forking
PIDFile=/run/hiawatha.pid
ExecStartPre=/usr/sbin/wigwam
ExecStart=/usr/sbin/hiawatha

[Install]
WantedBy=multi-user.target
