# disk-indicator/disk-indicator.service
[Unit]
Description=Disk Indicator

[Service]
ExecStart=/usr/bin/disk_indicator

[Install]
WantedBy=multi-user.target

