# mongodb2.4/mongodb.service
[Unit]
Description=High-performance, schema-free document-oriented database
After=network.target

[Service]
User=mongodb2_4
ExecStart=/usr/bin/mongod2.4 --quiet --config /etc/mongodb2.4.conf

[Install]
WantedBy=multi-user.target
