# udhcp/udhcpd.service
[Unit]
Description=udhcpd
Wants=network.target
Before=network.target

[Service]
ExecStart=/usr/sbin/udhcpd /etc/udhcpd.conf
ExecStop=/bin/kill -TERM $MAINPID
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
