# systemd-nexus-units/nexus.service
[Unit]
Description=A Maven Respository Server of Nexus
After=network.target

[Service]
Type=forking
EnvironmentFile=-/etc/conf.d/nexus
ExecStart=/opt/nexus/bin/nexus start
ExecStop=/opt/nexus/bin/nexus stop
ExecReload=/opt/nexus/bin/nexus restart
PIDFile=/run/nexus/nexus.pid
User=nexus

[Install]
WantedBy=multi-user.target

