# systemd-nexus-units/nexus.service
[Unit]
Description=A Maven Respository Server of Nexus
After=network.target

[Service]
Type=simple
ExecStart=/opt/nexus/bin/nexus start
ExecStop=/opt/nexus/bin/nexus stop
PIDFile=/var/run/nexus/nexus.pid

[Install]
WantedBy=multi-user.target

