# systemd-loop-swapfile/loop-swapfile.service
[Unit]
Description=Mount the swap file via a loop device
After=local-fs.target
Before=swap.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/sbin/losetup -f
ExecStart=/sbin/losetup /dev/loop0 /var/swapfile
ExecStart=/sbin/swapon /dev/loop0
ExecStop=/sbin/swapoff /dev/loop0
ExecStop=/sbin/losetup -d /dev/loop0

[Install]
Alias=local-fs.target.wants/loop-swapfile.service
