# newsd/newsd.service
[Unit]
Description=Standalone local NNTP server for local forum serving on a single server
After=network.target

[Service]
PrivateTmp=true
ExecStart=/usr/sbin/newsd -c /etc/newsd.conf -f

[Install]
WantedBy=multi-user.target
