# python2-pywo-svn/pywo.service
[Unit]
Description=pywo window organizer
Documentation=https://code.google.com/p/pywo/wiki/PywoConfiguration

[Service]
Type=simple
Environment=HOME=%h
Environment=DISPLAY=:0
Environment=XAUTHORITY=%h/.Xauthority
ExecStart=/usr/bin/pywo --daemon
KillMode=process

[Install]
WantedBy=multi-user.target

