# vmware_clock_sync/vmware_clock_sync.service
[Unit]
Description=Fix system clock drift that occurs when VMWare host sleeps.
Wants=vmware_clock_sync.timer

[Service]
Type=simple
Nice=19
ExecStart=/usr/bin/hwclock --hctosys --utc

[Install]
WantedBy=basic.target
