# plone-zeo/plone.service
[Unit]
Description=Plone web application server
After=network.target

[Service]
Type=forking
User=plone
ExecStart=/opt/Plone/zeocluster/bin/plonectl start
ExecStop=/opt/Plone/zeocluster/bin/plonectl stop

[Install]
WantedBy=multi-user.target