# ruby-eye/eye.service
[Unit]
Description=Eye

[Service]
Type=simple
PIDFile=/run/eye.pid
ExecStart=/usr/lib/ruby/gems/2.0.0/bin/eye -c /etc/eye/eye.conf -P /run/eye.pid -l /var/log/eye --log-level info -D
ExecStop=/usr/lib/ruby/gems/2.0.0/bin/eye terminate

[Install]
WantedBy=multi-user.target
