# youtrack-daemon/youtrack.service
[Unit]
Description=Issue Tracker and Agile Project Management Tool
After=network.target

[Service]
User=youtrack
Group=youtrack
Type=simple
EnvironmentFile=/etc/conf.d/youtrack
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=youtrack
ExecStart=/bin/sh -c 'eval $YOUTRACK_COMMAND_LINE'

[Install]
WantedBy=multi-user.target