# yaskkserv/yaskkserv.service
[Unit]
Description = Yet another SKK Japanese input method dictionary server

[Service]
Type = forking
PIDFile = /run/yaskkserv.pid
ExecStart = /usr/lib/systemd/scripts/yaskkserv-start

[Install]
WantedBy = multi-user.target
