# yaskkserv/yaskkserv.service
[Unit]
Description = Yet another SKK Japanese input method dictionary server

[Service]
Type = forking
User = nobody
EnvironmentFile = /etc/conf.d/yaskkserv
ExecStart = /usr/sbin/yaskkserv_hairy $YASKKSERV_ARGS $YASKKSERV_JISYO

[Install]
WantedBy = multi-user.target
