# dump1090_mr-git/dump1090.service
[Unit]
Description=RtlSdr mode S decoder
Wants=network-online.target
After=syslog.target network-online.target

[Service]
Restart=always
RestartSec=30
ExecStart=/usr/bin/dump1090 --quiet --net --net-ro-size 500 --net-ro-rate 5 --net-buffer 5
PIDFile="/run/dump1090.pid"

[Install]
WantedBy=multi-user.target
