# systemd-catalyst-pxp-mux-units/catalyst_pxp.service
[Unit]
Description=Deactivate the discrete graphic
After=systemd-modules-load.service
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/usr/bin/catalyst_pxp

[Install]
RequiredBy=multi-user.target
