# nzbdrone-torrents/nzbdrone.service
[Unit]
Description=NzbDrone-torrents
After=network.target

[Service]
User=nzbdrone
Group=nzbdrone
ExecStart=/usr/bin/nzbdrone -nobrowser -data=/var/lib/nzbdrone
Type=simple
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
