# urbackup-server/urbackup-server.service
[Unit]
Description=Backup Server
[Service]
ExecStart=/usr/sbin/start_urbackup_server --no_daemon
ExecStop="/usr/bin/killall -9 urbackup_srv"
[Install]
WantedBy=multi-user.target







