# urbackup-server/urbackup-server.service

[Unit]
Description=Backup Server
[Service]
ExecStart=/usr/sbin/start_urbackup_server --no_daemon
ExecStop="/usr/bin/killall -9 urbackup_srv"
[Install]
WantedBy=multi-user.target

# This service file runs an SSH daemon that forks for each incoming connection.
# If you prefer to spawn on-demand daemons, use sshd.socket and sshd@.service.










