# quassel-client-git/quassel-core.service
[Unit]
Description=Quassel core service
After=network.target

[Service]
Type=simple
EnvironmentFile=-/etc/conf.d/quassel-core
User=${user}
ExecStart=/usr/bin/quasselcore --configdir="${configdir}" --listen="${listen}"

[Install]
WantedBy=multi-user.target
