# qafooprofiler-daemon/qprofd.service
[Unit]
Description=Qafoo Profiler Daemon
After=network.target

[Service]
Type=simple
User=qprofd
Group=qprofd
ExecStart=/usr/bin/qprofd -pidfile=/run/qprofd/qprofd.pid
PIDFile=/run/qprofd/qprofd.pid
Restart=always

[Install]
WantedBy=multi-user.target
