# connman-dispatcher-chrony/connman-dispatcher.service
[Unit]
Description=Connman Dispatcher Daemon
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/connman-dispatcher -p /etc/connman-dispatcher.d/

[Install]
WantedBy=multi-user.target
