# graylog2-web-interface/graylog2-web-interface.service
[Unit]
Description=Graylog2 web-interface 
After=network.target

[Service]
User=nobody
EnvironmentFile=/etc/conf.d/graylog2-web-interface.conf
ExecStart=/usr/bin/java -Dpidfile.path=/dev/null -Duser.dir=${PATH} -Dhttp.port=${HTTP_PORT} -Dhttp.address=${HTTP_ADDRESS} -cp "${PATH}/lib/graylog2-web-interface.graylog2-web-interface-0.92.0.jar:${PATH}/lib/*" play.core.server.NettyServer

[Install]
WantedBy=multi-user.target
