# graylog2-web-interface/graylog2-web-interface.service
[Unit]
Description=Graylog2 web-interface 
After=network.target

[Service]
EnvironmentFile=/etc/conf.d/graylog2-web-interface.conf
PIDFile=/run/graylog2.pid
ExecStart=/usr/bin/java -Duser.dir=/usr/share/webapps/graylog2-web-interface/ -Dhttp.port=${HTTP_PORT} -Dhttp.address=${HTTP_ADDRESS} -cp "./lib/*" play.core.server.NettyServer

[Install]
WantedBy=multi-user.target
