# graylog2-web-interface/graylog2-web-interface.service
[Unit]
Description=Graylog2 web-interface 
After=network.target

[Service]
EnvironmentFile=/etc/conf.d/graylog2-web-interface
PIDFile=/run/graylog2.pid
ExecStart=/usr/bin/java -Duser.dir=/tmp -Dhttp.port=$HTTP_PORT -Dhttp.address=$HTTP_ADDRESS -cp "./lib/*" play.core.server.NettyServe

[Install]
WantedBy=multi-user.target
