# graylog2-web-interface/graylog2.cron
*/15 * * * * cd /usr/share/webapps/graylog2/ && bin/rake RAILS_ENV=production streamalarms:send > /dev/null