# psi-toolkit-git/psi-server.service
[Unit]
Description=PSI-Toolkit server
After=network.target

[Service]
ExecStart=/usr/bin/psi-server
User=psiserver
Group=psiserver

[Install]
WantedBy=multi-user.target