# gitorious-git/gitorious.service
[Unit]
Description=Gitorious
Wants=gitorious-resque.service git-daemon.socket

[Service]
User=gitorious
Group=gitorious
WorkingDirectory=/usr/share/webapps/gitorious
ExecStart=/usr/bin/bundle exec rails server -e production
Type=simple

[Install]
WantedBy=multi-user.target
