# gitorious-git/gitorious-sphinx.service
[Unit]
Description=Gitorious Sphinx
Wants=mysqld.service

[Service]
User=gitorious
Group=gitorious
ExecStart=/usr/bin/searchd -c /usr/share/webapps/gitorious/config/production.sphinx.conf
Type=forking
PidFile=/usr/share/webapps/gitorious/log/searchd.production.pid

[Install]
WantedBy=multi-user.target

