# gitorious-git/gitorious-resque.service
[Unit]
Description=Gitorious Resque Worker
Requires=redis.service

[Service]
User=gitorious
Group=gitorious
WorkingDirectory=/usr/share/webapps/gitorious
Environment=RAILS_ENV=production QUEUE=*
ExecStart=/usr/bin/bundle exec rake resque:work
Type=simple

[Install]
WantedBy=multi-user.target
