# systemd-metalog/metalog.service
[Unit]
Description=Alternative Syslog Daemon
Documentation=man:metalog(8)

[Service]
ExecStart=/usr/sbin/metalog
ExecReload=/bin/kill -HUP $MAINPID
