# obelisk-git/obbalancer.service
[Unit]
Description=Obelisk balancer
Requires=obworker.service
After=obworker.service

[Service]
Type=forking
PIDFile=/run/obbalancer.pid
User=ob
Group=ob
ExecStart=/usr/bin/obbalancer
Restart=always

[Install]
WantedBy=multi-user.target
