# obelisk-git/obelisk.service
[Unit]
Description=Obelisk worker
Requires=sx-initchain.service
After=sx-initchain.service
OnFailureJobMode=ignore-requirements

[Service]
Type=forking
PIDFile=/run/obworker.pid
User=ob
Group=ob
ExecStart=/usr/bin/obworker
Restart=always

[Install]
WantedBy=multi-user.target
