# obelisk-git/obelisk.service
[Unit]
Description=Obelisk worker
Requires=sx-initchain.service
After=sx-initchain.service
OnFailureJobMode=ignore-requirements

[Service]
Type=forking
PIDFile=/run/obworker.pid
User=obelisk
Group=obelisk
ExecStart=/usr/bin/obworker
Restart=always
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
