# mlocate-systemd/mlocate-update.service
[Unit]
Description=Update mlocate database

[Service]
Nice=19
IOSchedulingClass=2
IOSchedulingPriority=7
ExecStart=/usr/bin/updatedb

[Install]
WantedBy=cron-daily.target
