# autosshd/autosshd.service
[Unit]
Description=Automatically restart SSH sessions and tunnels
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/autossh #<SSH arguments>

[Install]
WantedBy=multi-user.target
