# magic-motd/magic-motd-ip.service
[Unit]
Description=Message of the day generator (IP Refresher)
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/magic-motd-ip

[Install]
WantedBy=multi-user.target
