# qtsixa/sixad.service
[Unit]
Description=SixAxis daemon
After=bluetooth.target

[Service]
Type=simple
ExecStart=/usr/bin/sixad --start
ExecStop=/usr/bin/sixad --stop

[Install]
WantedBy=multi-user.target
