# selinux-cronie/crontab
### /etc/crontab: Master crontab for system `crond`

PATH=/sbin:/bin:/usr/sbin:/usr/bin
SHELL=/bin/bash
HOME=/var/empty

# Format: Minute Hour Date Month Wday User Command

01  *  *  *  *  root  run-cron /etc/cron.hourly
02  4  *  *  *  root  run-cron /etc/cron.daily
22  4  *  *  0  root  run-cron /etc/cron.weekly
42  4  1  *  *  root  run-cron /etc/cron.monthly

### /etc/crontab: Master crontab for system `crond`
