# antivir/antivir.service
[Unit]
Description=AntiVir Guard Daemon

[Service]
Type=forking
ExecStart=/usr/sbin/avguard start
ExecStop=/usr/sbin/avguard stop

[Install]
WantedBy=multi-user.target
