# tlp/tlp-sleep.service
[Unit]
Description=TLP resume reinstate
Befor=sleep.target
StopWhenUnneeded=yes

[Service]
User=%I
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/tlp start

[Install]
WantedBy=sleep.target
