# xfluxd/xfluxd.service
[Unit]
Description=Daemon for xflux
After=display-manager.service

[Service]
Type=forking
EnvironmentFile=/etc/xfluxd.conf
ExecStart=/usr/bin/xflux -l $LAT -g $LON -k $TEMP
Restart=always

[Install]
WantedBy=graphical.target
