# astrapiremoted/astrapiremoted.service
[Unit]
Description=Astrapi Remote Daemon

[Service]
ExecStart=/opt/astrapi/bin/astrapiremoted

[Install]
WantedBy=multi-user.target
