# mediacrush-server-git/mediacrush.service
[Unit]
Description=MediaCrush gunicorn daemon
Wants=mediacrush-daemon.service redis.service

[Service]
User=http
WorkingDirectory=/var/lib/mediacrush
ExecStart=/usr/bin/gunicorn app:app -b 127.0.0.1:35441

[Install]
WantedBy=multi-user.target
