# mediacrush-server-git/mediacrush-daemon.service
[Unit]
Description=MediaCrush ffmpeg daemon

[Service]
User=http
WorkingDirectory=/var/lib/mediacrush
ExecStart=/usr/bin/python2 daemon.py

[Install]
WantedBy=multi-user.target
