# btsync/btsync.service
[Unit]
Description=BTSync
After=network.target

[Service]
User=root
ExecStart=/usr/share/btsync/btsync --config /etc/btsync/btsync.json --nodaemon

[Install]
WantedBy=multi-user.target
