# btsync/btsync@.service
[Unit]
Description=BTSync for %i
After=network.target
 
[Service]
User=%i
ExecStart=/usr/bin/btsync --config /home/%i/.btsync/settings.json --nodaemon
 
[Install]
WantedBy=multi-user.target
