# batteryd-git/batteryd.service
[Unit]
Description=Batteryd - the smallest battery daemon ever
After=syslog.target 

[Service]
User=user
Environment=DISPLAY=:0
Type=simple
ExecStart=/usr/bin/batteryd 

[Install]
WantedBy=multi-user.target

