# autobacklight-git/autobacklight.service
[Unit]
Description=Autobacklight Daemon

[Service]
ExecStart=/usr/bin/autobacklight
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=always

[Install]
WantedBy=multi-user.target
