# hptwebgui/hptwebgui.service
[Unit]
Description = HighPoint RocketRAID webGUI daemon

[Service]
Type = forking
ExecStart = /usr/bin/hptsvr /etc/hptcfg
ExecStop = /usr/bin/hptsvr --kill /etc/hptcfg
Restart = on-abort

[Install]
WantedBy=multi-user.target
