# hostsblock/hostsblock.timer
[Unit]
Description=hostsblock.timer: Daily Timer for hostsblock.service

[Timer]
OnBootSec=15m
OnCalendar=daily
Persistent=true

[Install]
WantedBy=multi-user.target
