# hostsblock/hostsblock.service
[Unit]
Description=Block bad hosts
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/hostsblock
