# gdm-old/gdm.service
[Unit]
Description=Gnome Display Manager
Conflicts=getty@tty1.service
After=systemd-user-sessions.service getty@tty1.service

[Service]
ExecStart=/usr/sbin/gdm --nodaemon
ExecReload=/usr/sbin/gdm-restart
ExecStop=/usr/sbin/gdm-stop
RestartSec=1s
TimeoutStopSec=5s

[Install]
Alias=display-manager.service
