# gdm-old/gdm.service
[Unit]
Description=Gnome Display Manager
After=systemd-user-sessions.service

[Service]
ExecStart=/usr/bin/gdm --nodaemon
ExecReload=/usr/bin/gdm-safe-restart
ExecStop=/usr/bin/gdm-stop
IgnoreSIGPIPE=no
RestartSec=1s
TimeoutStopSec=5s

[Install]
Alias=display-manager.service
