# cpufreq-systemd/cpufreq.service
[Unit]
Description=CPU frequency scaling via cpufrequtils

[Service]
ExecStart=/usr/bin/cpufreq-systemd.sh
ExecReload=/usr/bin/cpufreq-systemd.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
