# cudaminer/cudaminer.service
[Unit]
Description=CUDA accelerated mining application
Documentation=file:///usr/share/docs/cudaminer/README
Documentation=https://bitcointalk.org/index.php?topic=167229.0
Documentation=https://docs.google.com/spreadsheet/ccc?key=0Aj3vcsuY-JFNdHR4ZUN5alozQUNvU1pyd2NGeTNicGc
ConditionFileNotEmpty=/etc/cudaminer.conf
After=network.target mining-proxy.service syslog.target

[Service]
ExecStart=/usr/bin/cudaminer --config /etc/cudaminer.conf --syslog

[Install]
WantedBy=multi-user.target
