# puppetdb/puppetdb.service
[Unit]
Description=puppetdb Daemon
After=syslog.target network.target

[Service]
Type=simple
User=puppetdb
PIDFile=/var/run/puppetdb.pid
ExecStart=/usr/bin/java \
          -Xmx192m -XX:+HeapDumpOnOutOfMemoryError -XX:HeapDumpPath=/var/log/puppetdb/puppetdb-oom.hprof  \
          -jar /usr/share/puppetdb/puppetdb.jar \
          services -c /etc/puppetdb/conf.d $@
  
ExecStop=/bin/kill $MAINPID

[Install]
WantedBy=multi-user.target
