# hqplayer-embedded/hqplayerd@.service
[Unit]
Description=HQ player embedded
After=network.target sound.target

[Service]
ExecStart=/usr/bin/hqplayerd --daemon
Restart=always
User=%I

[Install]
WantedBy=multi-user.target

