# cdemu-systemd/cdemu.service
[Unit]
Description=CD/DVD-ROM device emulator daemon

[Service]
EnvironmentFile=/etc/conf.d/cdemud
ExecStartPre=/sbin/modprobe vhba
ExecStart=/usr/bin/cdemud -c $CTL_DEVICE -b $BUS -n $NUM_DEVICES -a $AUDIO_DRIVER -l $LOG_FILE
ExecStop=/bin/kill -HUP $MAINPID
ExecStopPost=/usr/bin/rmmod -f vhba

[Install]
WantedBy=multi-user.target
