# rygel-systemd/rygel@.service
[Unit]
Description=Rygel DLNA server
After=syslog.target local-fs.target network.target

[Service]
ExecStart=/usr/bin/rygel
Restart=always
User=%I

[Install]
WantedBy=multi-user.target

