# bliss/bliss.service
[Unit]
Description=Bliss Daemon

[Service]
Type=simple
User=bliss
Group=bliss
ExecStart=/opt/bliss/bin/bliss.sh

[Install]
WantedBy=multi-user.target
