# eclim-git/systemd_eclimd.service
[Unit]
Description=Eclim
After=graphical.target

[Service]
Environment=DISPLAY=:0
ExecStart=/usr/share/eclipse/eclimd
