# synergy-1.3/synergys-1.3.service
[Unit]
Description=Synergy v1.3 Server Daemon
After=network.target

[Service]
Type=forking
ExecStart=/usr/bin/synergys-1.3 --config /etc/synergy-1.3.conf

[Install]
WantedBy=multi-user.target
