# irqbalance-novena/irqbalance.service
[Unit]
Description=CPU Interrupt Request Balancer
After=sysinit.target

[Service]
ExecStart=/usr/bin/irqbalance
Type=forking

[Install]
WantedBy=multi-user.target
