# grafana/grafana.service
[Unit]
Description=Grafana service
After=network.target                                     

[Service]
User=grafana
ExecStart=/usr/bin/grafana --config /etc/grafana/grafana.ini web

[Install]
WantedBy=multi-user.target
