# zabbix-server/zabbix-server.service
[Unit]
Description=Zabbix server
After=syslog.target network.target postgresql.service

[Service]
Type=oneshot
ExecStart=/usr/sbin/zabbix_server
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target

