# opensmtpd-snapshot/smtpd.service
[Unit]
Description=A FREE implementation of the server-side SMTP protocol
After=syslog.target
Wants=network.target
Before=network.target

[Service]
ExecStart=/usr/sbin/smtpd -d

[Install]
WantedBy=multi-user.target
