# geoip-update-systemd/geoip-update.timer
[Timer]
OnUnitInactiveSec=1w
OnStartupSec=1h
AccuracySec=1h
