# geoip-update-systemd/geoip-update.service
[Service]
Type=oneshot
WorkingDirectory=/usr/share/GeoIP

ExecStart=/bin/bash -c 'wget -nv http://geolite.maxmind.com/download/geoip/database/GeoLiteCountry/GeoIP.dat.gz -O - | gunzip -c - > GeoIP.dat'
ExecStart=/bin/bash -c 'wget -nv http://geolite.maxmind.com/download/geoip/database/GeoLiteCity.dat.gz -O - | gunzip -c - > GeoLiteCity.dat'
ExecStart=/bin/bash -c 'wget -nv http://geolite.maxmind.com/download/geoip/database/GeoIPv6.dat.gz -O - | gunzip -c > GeoIPv6.dat'
ExecStart=/bin/bash -c 'wget -nv http://geolite.maxmind.com/download/geoip/database/GeoLiteCityv6-beta/GeoLiteCityv6.dat.gz -O - | gunzip -c > GeoLiteCityv6.dat'
