# wol-systemd/wol@.service
[Unit]
Description=Wake-on-LAN for %i
Requires=network.target
After=sys.subsystem-net-devices-%i.device

[Service]
ExecStart=/usr/sbin/enable-wol %i
Type=oneshot

[Install]
WantedBy=multi-user.target
