# nagios/nagios.service
[Unit]
Description=Nagios Monitoring System
After=network.target

[Service]
ExecStart=/usr/bin/nagios /etc/nagios/nagios.cfg

[Install]
WantedBy=multi-user.target
