# i8kutils-bzr/i8kmon.service
[Unit]
Description=Dell laptop thermal monitoring
ConditionPathExists=/proc/i8k
After=systemd-modules-load.service

[Service]
ExecStart=/usr/bin/i8kmon -d
RestartSec=4
Restart=always

[Unit]
DefaultDependencies=no

[Install]
WantedBy=sysinit.target
