# core/core-emulator.service
[Unit]
Description=Common Open Research Emulator Service
After=network.target

[Service]
Type=forking
PIDFile=/run/cored.pid
ExecStart=/usr/bin/python2 /usr/sbin/cored.py -d -i /run/cored.pid

[Install]
WantedBy=multi-user.target
