# shadowsocks-nodejs/shadowsocks-nodejs.service
[Unit]
Description=A lightweight tunnel proxy.

[Service]
ExecStart=/usr/bin/sslocal -c /usr/lib/node_modules/shadowsocks/config.json
StandardOutput=null

[Install]
WantedBy=multi-user.target
