# mylar-git/mylar-user.service
[Unit]
Description=Mylar Daemon

[Service]
ExecStart=/usr/bin/env python2 /opt/mylar/Mylar.py --config %h/.mylar/config.ini --datadir %h/.mylar

[Install]
WantedBy=default.target
