# mylar-git/mylar-system.service
[Unit]
Description=Mylar Daemon
After=network.target

[Service]
User=mylar
Group=mylar
ExecStart=/usr/bin/env python2 /opt/mylar/Mylar.py --config /opt/mylar/config.ini --datadir /opt/mylar

[Install]
WantedBy=multi-user.target
