# brewpi-git/brewpi.cron
stderrpath="/home/brewpi/logs/stderr.txt"
stdoutpath="/home/brewpi/logs/stdout.txt"
scriptpath="/home/brewpi"
entries="brewpi ~wifichecker"
# entry:brewpi
* * * * * brewpi python2 $scriptpath/brewpi.py --checkstartuponly --dontrunfile $scriptpath/brewpi.py 1>/dev/null 2>>$stderrpath; [ $? != 0 ] && python2 -u $scriptpath/brewpi.py 1>$stdoutpath 2>>$stderrpath &

