# zram/zdisk@.service
[Unit]
Description=zRAM disk %i

[Service]
Type=oneshot
RemainAfterExit=true

ExecStart=/usr/sbin/zdisk.sh start %i
ExecStop=/usr/sbin/zdisk.sh stop %i

[Install]
WantedBy=local-fs.target
