# systemd-mongodb-units/mongodb.service
[Unit]
Description=MongoDB Server
After=network.target

[Service]
User=mongodb
ExecStart=/usr/bin/mongod --config /etc/mongodb.conf

[Install]
WantedBy=multi-user.target
