# autorestartnm/timer-quarterhourly.target
[Unit]
Description=Quarter-Hourly Timer Target
StopWhenUnneeded=yes
