# autorestartnm/timer-quarterhourly.timer
[Unit]
Description=Quarter-Hourly Timer

[Timer]
OnBootSec=8min
OnUnitActiveSec=15min
Unit=timer-quarterhourly.target

[Install]
WantedBy=basic.target
