# fdnsc/fdnsc.service
[Unit]
Description=freedns.afraid.org dynamic dns update client
After=network.target

[Service]
ExecStart=/usr/sbin/fdnsc
ExecStop=/usr/sbin/fdnsc stop

[Install]
WantedBy=multi-user.target
