# openntpd-fast/openntpd-fast.service
[Unit]
Description=OpenNTP Daemon
After=network.target
Conflicts=openntpd

[Service]
ExecStart=/usr/bin/ntpd -ds

[Install]
WantedBy=multi-user.target
