# jboss-as/jboss-as.service
[Unit]
Description=JBoss Application Server

[Service]
Type=simple
User=jboss-as
PIDFile=/run/jboss-as/jboss-as.pid
ExecStart=/usr/share/jboss-as/bin/systemd-wrapper.sh
ExecReload=/usr/share/jboss-as/bin/jboss-cli.sh --connect command=:reload
ExecStop=/usr/share/jboss-as/bin/jboss-cli.sh --connect command=:shutdown

[Install]
WantedBy=multi-user.target
