# ruby-redmine2/redmine2.service
[Unit]
Description=Redmine webserver
After=network.target

[Service]
User=redmine
ExecStart=/usr/bin/ruby /var/lib/redmine2/script/rails server webrick -e production
ExecStop=/bin/kill -INT $MAINPID

[Install]
WantedBy=multi-user.target