# gitlab/gitlab-sidekiq.service
[Unit]
Description=Gitlab Sidekiq
Wants=redis.service

[Service]
Type=forking
User=gitlab
WorkingDirectory=/usr/share/webapps/gitlab
Environment=RAILS_ENV=production

ExecStart=/usr/bin/bundle exec rake sidekiq:start
ExecStop=/usr/bin/bundle exec rake sidekiq:stop

[Install]
WantedBy=multi-user.target
