# gitlab/gitlab.service
[Unit]
Description=Gitlab Rails server
After=network.target mysql.service

[Service]
Type=simple
User=gitlab
PIDFile=/run/gitlab.pid
WorkingDirectory=/usr/share/webapps/gitlab
Environment=RAILS_ENV=production

ExecStart=/usr/bin/bundle exec rails server

[Install]
WantedBy=multi-user.target
