# foldingathome-gpu/foldingathome-gpu@.service
[Unit]
Description=Folding@Home GPU client (GPU %i)

[Service]
ExecStart=/usr/bin/foldingathome-gpu-wrapper -c /etc/foldingathome-gpu/client.cfg -w /var/lib/foldingathome-gpu/gpu%i -g %i

[Install]
WantedBy=multi-user.target
