# foldingathome-gpu/foldingathome-gpu@.service
[Unit]
Description=Folding@Home GPU client (GPU %i)

[Service]
ExecStart=/usr/bin/foldingathome-gpu-wrapper %i

[Install]
WantedBy=multi-user.target
