# systemd-loop-swapfile-auto/systemd-loop-swapfile.service
[Unit]
Description=Mount the swap file via a loop device
After=local-fs.target
Before=swap.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/lib/systemd/scripts/systemd-loop-swapfile.sh
ExecStop=/usr/bin/swapoff $(/usr/bin/cat /tmp/swaploopdev)

[Install]
Alias=local-fs.target.wants/systemd-loop-swapfile.service
