# hlds/hlds.service
[Unit]                                                                                                                         
Description=Half-Life Dedicated Server
After=network.target

[Service]
Type=Forking
EnvironmentFile=/etc/conf.d/hlds
User=hlds
WorkingDirectory=/opt/hlds
ExecStart=/opt/hlds/hlds_run ${params}
ExecStop=/bin/kill $MAINPID
Restart=always

[Install]
WantedBy=multi-user.target

