# nfs-ganesha/nfs-ganesha.service
[Unit]
Description=User space NFS server
Documentation=https://github.com/nfs-ganesha/nfs-ganesha/wiki

[Service]
ExecStart=/usr/bin/ganesha.nfsd
Type=forking
PIDFile=/var/run/ganesha.pid

[Install]
WantedBy=multi-user.target
