# cisco-vpnclient-amdphenom/ciscovpn@.service
[Unit]
Description=Cisco VPN Client
After=network.target

[Service]
ExecStart=/usr/bin/vpnclient connect "%I"
ExecStop=/usr/bin/vpnclient disconnect

[Install]
WantedBy=multi-user.target
