# ffdecsawrapper-git-digital-devices/ffdecsawrapper.service
[Unit]
Description=FFdecsa empowered softcam
After=oscam.service 
Wants=oscam.service

[Service]
Type=simple
User=root
Nice=-2
EnvironmentFile=/etc/conf.d/ffdecsawrapper
ExecStartPre=/usr/bin/modprobe dvbloopback $DVBLOOPBACK_ARGS
TimeoutSec=2
ExecStart=/usr/bin/ffdecsawrapper $FFDECSAWRAPPER_ARGS --cam-dir $CAMDIR --log $LOGFILE

[Install]
WantedBy=multi-user.target
