# vdradmin-am/vdradmind.service
[Unit]
Description=vdradmin-am VDR web interface
After=vdr.service

[Service]
Type=forking
PIDFile=/run/vdradmind.pid
ExecStart=/usr/bin/vdradmind --pid /run/vdradmind.pid
ExecStop=-/usr/bin/vdradmind --kill

[Install]
WantedBy=multi-user.target
