# zfs-utils/zfs-utils.service
[Unit]
Description=Imports and mounts ZFS pools.
Documentation=man:zfs(8) man:zpool(8)
DefaultDependencies=no
Conflicts=shutdown.target umount.target

[Service]
Type=oneshot
ExecStart=/usr/sbin/zfs mount -a
ExecStop=/usr/sbin/zfs unmount -a
RemainAfterExit=yes

[Install]
WantedBy=sysinit.target
