# plymouth-systemd-git/gdm-plymouth.service
[Unit]
Description=GNOME Display Manager
After=systemd-user-sessions.service

Conflicts=getty@tty1.service plymouth-quit.service
After=getty@tty1.service plymouth-quit.service

[Service]
ExecStart=/usr/sbin/gdm -nodaemon
StandardOutput=syslog

[Install]
WantedBy=graphical.target
