# shinken/shinken-scheduler.service
[Unit]
Description=Shinken Scheduler 
After=syslog.target

[Service]
Type=forking
ExecStart=/usr/bin/shinken-scheduler -dc /etc/shinken/schedulerd.ini

[Install]
WantedBy=multi-user.target
