# ceilometer/ceilometer-collector.service
[Unit]
Description=Ceilometer Collector service

[Service]
ExecStart=/usr/bin/ceilometer-collector
User=ceilometer

[Install]
WantedBy=multi-user.target
