# ceilometer/ceilometer-alarm-notifier.service
[Unit]
Description=Ceilometer Alarm Notifier

[Service]
ExecStart=/usr/bin/ceilometer-alarm-notifier --log-dir=/var/log/ceilometer
User=ceilometer

[Install]
WantedBy=multi-user.target
