# ceilometer/ceilometer-api.service
[Unit]
Description=Ceilometer API service

[Service]
ExecStart=/usr/bin/ceilometer-api --log-dir=/var/log/ceilometer
User=ceilometer

[Install]
WantedBy=multi-user.target
