# ceilometer/ceilometer-agent-notification.service
[Unit]
Description=Ceilometer Agent Notification

[Service]
ExecStart=/usr/bin/ceilometer-agent-notification
User=ceilometer

[Install]
WantedBy=multi-user.target
