# ceilometer/ceilometer-agent-central.service
[Unit]
Description=Ceilometer Agent Central

[Service]
ExecStart=/usr/bin/ceilometer-agent-central
User=ceilometer

[Install]
WantedBy=multi-user.target
