# ceilometer/ceilometer-alarm-evaluator.service
[Unit]
Description=Ceilometer Alarm Evaluator

[Service]
ExecStart=/usr/bin/ceilometer-alarm-evaluator --log-dir=/var/log/ceilometer
User=ceilometer

[Install]
WantedBy=multi-user.target
