# ceilometer/ceilometer-agent-compute.service
[Unit]
Description=Ceilometer Agent Compute

[Service]
ExecStart=/usr/bin/ceilometer-agent-compute
User=ceilometer

[Install]
WantedBy=multi-user.target
