# wildcard_dns_proxy/wildcard_dns_proxy.service
[Unit]
Description=Wildcard DNS proxy
After=rc-local.service

[Service]
Type=simple
ExecStart=/usr/bin/wildcard_dns_proxy -s 8.8.8.8
SysVStartPriority=99

[Install]
WantedBy=multi-user.target
