# mysql-wsrep/mysqld-bootstrap.service
[Unit]
Description=wsrep-enabled MySQL server(bootstrap)
After=syslog.target
Conflicts=mysqld.service

[Service]
User=mysql
Group=mysql

ExecStart=/usr/bin/mysqld --pid-file=/run/mysqld/mysqld.pid --wsrep_new_cluster
ExecStartPost=/usr/bin/mysqld-post

Restart=always
PrivateTmp=true

[Install]
WantedBy=multi-user.target
