# systemd-hddapm/hddapm.service
[Unit]
Description=Sets the Advanced Power Management level for /dev/sda

[Service]
Type=oneshot
ExecStart=/usr/bin/hddapm post
