# potcoin-daemon/potcoind@.service
[Unit]
Description=Potcoin Daemon
After=network.target 

[Service]
User=%i
Group=users
Type=simple
ExecStart=/usr/bin/potcoind

[Install]
WantedBy=multi-user.target
