# ruby-redmine/redmine.service
[Unit]
Description=Redmine webserver
After=network.target

[Service]
User=redmine
ExecStart=/usr/bin/ruby-1.8 /var/lib/redmine/script/server webrick -e production
ExecStop=/bin/kill -INT $MAINPID

[Install]
WantedBy=multi-user.target
