# quasselcore/quassel.service
[Unit]
Description=Quassel core service
After=network.target

[Service]
Type=simple
EnvironmentFile=-/etc/conf.d/quassel
User=quassel
ExecStart=/usr/bin/quasselcore --configdir=/var/lib/quassel --listen=${LISTEN}

[Install]
WantedBy=multi-user.target
