# raspi-off-button/raspi-off-button.service
[Unit]
Description=Raspberry PI Power-Off

[Service]
ExecStart=/usr/bin/raspi-off-button

[Install]
WantedBy=multi-user.target
