# pfix-srsd-git/pfix-srsd.service
[Unit]
Description=postfix srsd lookup service
 
[Service]
Type=simple
#User=pfix-srsd
#Group=pfix-srsd
#PIDFile=/run/pfix-srsd.pid
EnvironmentFile=-/etc/default/pfix-srsd
ExecStart=/usr/bin/pfix-srsd -U pfix-srsd -G pfix-srsd -f -v ${DOMAIN} ${SECRET_FILE}
 
[Install]
WantedBy=multi-user.target
