# icecc/icecc-scheduler.service
[Unit]
Description=icecc scheduler
Wants=network.target

[Service]
Type=forking
ExecStart=/usr/sbin/icecc-scheduler -d $ICECC_SCHEDULER_OPTIONS
EnvironmentFile=/etc/icecc.conf

[Install]
Alias=multi-user.target.wants/icecc-scheduler.service
