# icecc/icecc.service
[Unit]
Description=icecc worker
Wants=network.target

[Service]
Type=forking
ExecStart=/usr/sbin/iceccd -d $ICECCD_OPTIONS
EnvironmentFile=/etc/icecc.conf

[Install]
Alias=multi-user.target.wants/iceccd.service
