# icecc/icecc-scheduler.service
[Unit]
Description=icecc scheduler
Wants=network.target

[Service]
Type=forking
ExecStart=/usr/lib/icecc/icecc-scheduler

[Install]
Alias=multi-user.target.wants/icecc-scheduler.service
