# icecc/icecc-scheduler.service
[Unit]
Description=icecc scheduler
Wants=network.target
EnvironmentFile=/etc/icecc.conf

[Service]
Type=forking
ExecStart=/usr/lib/icecc/icecc-scheduler -d $ICECC_SCHEDULER_OPTIONS

[Install]
Alias=multi-user.target.wants/icecc-scheduler.service
