# icecc/icecc.service
[Unit]
Description=icecc worker
Wants=network.target

[Service]
Type=forking
ExecStart=/usr/lib/icecc/iceccd

[Install]
Alias=multi-user.target.wants/iceccd.service
