# icecc/icecc.service
[Unit]
Description=icecc worker
Wants=network.target
EnvironmentFile=/etc/icecc.conf

[Service]
Type=forking
ExecStart=/usr/lib/icecc/iceccd -d $ICECCD_OPTIONS

[Install]
Alias=multi-user.target.wants/iceccd.service
