# mooncoin-daemon/mooncoind@.service
[Unit]
Description=DogeCoin Daemon
After=network.target 

[Service]
User=%i
Group=users
Type=simple
ExecStart=/usr/bin/dogecoind

[Install]
WantedBy=multi-user.target
