# lastfmsubmitd/lastfmsubmitd.service
[Unit]
Description=LastFM Submit daemon
After=network.target

[Service]
ExecStart=/usr/bin/lastfmsubmitd --no-daemon

[Install]
WantedBy=multi-user.target
