# lastfmsubmitd/lastmp.service
[Unit]
Description=LastMP
After=mpd.service

[Service]
ExecStart=/usr/bin/lastmp --no-daemon

[Install]
WantedBy=multi-user.target
