# chromiumos-bsp-daisy/thermal.service
[Unit]
Description=Chromium OS thermal monitoring daemon

[Service]
Type=simple
ExecStart=/usr/sbin/thermal.sh
RestartSec=10s
Restart=always

[Install]
WantedBy=sysinit.target
