# webdis-git/webdis.service
[Unit]
Description=Redis HTTP interface with JSON output
Requires=redis.service
After=redis.service

[Service]
Type=forking
ExecStart=/usr/bin/webdis /etc/webdis.prod.json
Restart=always
