# ekeyd/egd-linux.service
[Unit]
Description=EGD Enropy Daemon
After=network.target
# After=network.target,ekeyd.service

[Service]
ExecStart=/usr/bin/egd-linux -H 127.0.0.1 -p 8888 -D/run/egd-linux.pid
ExecPost=/usr/bin/sysctl kernel.random.write_wakeup_threshold=1024
PIDFile=/run/egd-linux.pid
KillMode=process

[Install]
WantedBy=multi-user.target

