# ekeyd/ekeyd.service
[Unit]
Description=Entropy Key Daemon
After=network.target

[Service]
ExecStart=/usr/bin/ekeyd
PIDFile=/run/ekeyd.pid
KillMode=process

[Install]
WantedBy=multi-user.target

