# freeipa/ipa-kpasswd.service
[Unit]
Description=FreeIPA Kpasswd Daemon

[Service]
ExecStart=/usr/sbin/ipa_kpasswd
Type=forking

[Install]
WantedBy=multi-user.target
