# ery-twinkle-broadcast-server-git/twinkle-broadcast-server.service
[Unit]
Description=Twinkle broadcast server

[Service]
ExecStart=/usr/bin/twinkle-broadcast-server

[Install]
WantedBy=multi-user.target
