# hybrid-video-ati-intel-systemd/ati-intel-hybrid.service
[Unit]
Description=ati-intel hybrid

[Service]
Type=oneshot
User=root
ExecStart=/usr/sbin/ati-intel-hybrid
StandardOutput=syslog

[Install]
WantedBy=multi-user.target
Requires=sys-kernel-debug.mount
