# kvm-opensuse/65-kvm.rules
################################################################################
# UDEV RULES FOR KVM
################################################################################

# Give the /dev/kvm to the kvm group
# Add your user account to that group to use KVM.
KERNEL=="kvm", NAME="%k", GROUP="kvm", MODE="0660"

# Do the same for network devices.
KERNEL=="tun", GROUP="kvm", MODE="0660"

# when the user loads the kvm module, automatically load kvm-intel or
# kvm-amd depending on the hardware
ACTION=="add", DEVPATH=="/module/kvm", RUN+="/bin/sh -c 'grep -q vmx /proc/cpuinfo   +&& /sbin/modprobe kvm-intel; grep -q svm /proc/cpuinfo && /sbin/modprobe kvm-amd'"
