# dudle-bzr/dudle.service
[Unit]
Description=Etherpad-Lite server
After=syslog.target
After=network.target

[Service]
Type=simple
User=http
Group=http
ExecStart=/usr/bin/etherpad-lite

[Install]
WantedBy=multi-user.target
