# cras-daisy/cras-state.service
[Unit]
Description=ChromeOS audio server state management
Requires=cras.service
After=cras.service

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/bin/crasctl init
ExecStart=/usr/bin/crasctl restore
ExecStop=/usr/bin/crasctl store
