# cras-daisy/cras.service
[Unit]
Description=ChromeOS audio server
Requires=sound.target dbus.service cras.socket cras-state.service
After=sound.target dbus.service cras.socket
Before=cras-state.service

[Service]
Type=dbus
BusName=org.chromium.cras
ExecStartPre=/usr/bin/chmod -R 777 /var/run/cras
ExecStart=/usr/bin/sg audio /usr/bin/cras
