# verlihub-git/verlihub.service
[Unit]
Description=Verlihub DC Server
After=network.target

[Service]
Type=forking
PIDFile=/etc/verlihub/pid

ExecStart=/usr/bin/vh --run
ExecStop=/usr/bin/vh --stop
ExecReload=/usr/bin/vh --restart

PrivateTmp=true
Restart=always

[Install]
WantedBy=multi-user.target 
