# freelan/freelan.service
[Unit]
Description=Open-source, multi-platform, peer-to-peer VPN software

[Service]
Type=forking
ExecStart=/usr/bin/freelan
Environment="HOME=/"

[Install]
WantedBy=multi-user.target
