# elasticsearch/elasticsearch.service
[Unit]
Description=ElasticSearch
After=network.target

[Service]
User=elasticsearch
ExecStart=/usr/share/elasticsearch/bin/elasticsearch -f -Des.path.conf=/etc/elasticsearch

[Install]
WantedBy=multi-user.target
