# elasticsearch/elasticsearch.service
[Unit]
Description=ElasticSearch FullTextIndex Engine
After=syslog.target

[Service]
Type=forking
ExecStart=/usr/share/elasticsearch/bin/service/elasticsearch start
ExecStop=/usr/share/elasticsearch/bin/service/elasticsearch stop

[Install]
WantedBy=multi-user.target
