# elasticsearch/elasticsearch.service
[Unit]
Description=ElasticSearch FullTextIndex Engine
After=syslog.target

[Service]
Type=forking
ExecStart=/opt/elasticsearch/bin/service/elasticsearch start
ExecStop=/opt/elasticsearch/bin/service/elasticsearch stop

[Install]
WantedBy=multi-user.target
