# csprng/csprngd.service
[Unit]
Description=Cryptographically Secure PseudoRandom Number Generator Daemon

[Service]
ExecStart=/usr/bin/csprngd --foreground

[Install]
WantedBy=multi-user.target
