# rstudio-server-git/rstudio-server.service
[Unit]
Description=RStudio Server daemon
After=syslog.target network.target

[Service]
Type=forking
PIDFile=/run/rstudio-server.pid
ExecStart=/usr/sbin/rstudio-server --server-user=rstudio --www-local-path=/srv/rstudio --server-daemonize
ExecStop=/bin/rm -f /run/rstudio-server.pid

[Install]
WantedBy=multi-user.target
