# rstudio-server-git/rstudio-server.service
[Unit]
Description=RStudio Server daemon
After=syslog.target network.target

[Service]

ExecStart=/usr/bin/rserver --server-daemonize=0
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=always

[Install]
WantedBy=multi-user.target
