# rstudio-server-git/rstudio-server.service
[Unit]
Description=RStudio Server
After=syslog.target network.target

[Service]
Type=forking
EnvironmentFile=/etc/conf.d/rstudio-server
ExecStartPre=/usr/lib/rstudio-server/bin/rserver --test-config
ExecStart=/usr/lib/rstudio-server/bin/rserver $RSTUDIO_SERVER_ARGS
ExecStop=/bin/kill -15 $MAINPID

[Install]
WantedBy=multi-user.target
