# smartladder/smartladder.service
[Unit]
Description=Smartladder Service
After=network.target

[Service]
Type=simple
ExecStart=/opt/smartladder/goagent+/proxy.py
ExecStop=/bin/kill -15 $MAINPID
StandardOutput=null

[Install]
WantedBy=multi-user.target

