# ptpd/ptpd.service
[Unit]
Description=Precision Time Protocol version 1 daemon

[Service]
ExecStart=/usr/bin/ptpd

[Install]
WantedBy=multi-user.target
