# plone/plone-standalone.service
[Unit]
Description=Plone
After=network.target
     
[Service]
Type=forking
ExecStart=/opt/plone/zinstance/bin/plonectl start
ExecStop=/opt/plone/zinstance/bin/plonectl stop
ExecReload=/opt/plone/zinstance/bin/plonectl restart
   
[Install]
WantedBy=multi-user.target
