# slurm-llnl/slurmctld.service
[Unit]
Description=Start/Stop the SLURMCTLD service
After=network.target auditd.service munged.service slurmdbd.service
Requires=munged.service

[Service]
Type=forking
ExecStart=/usr/bin/slurmctld
ExecStop=/bin/kill $MAINPID
ExecReload=/bin/kill -HUP $MAINPID
PIDFile=/run/slurmctld.pid

[Install]
WantedBy=multi-user.target
