# pure-ftpd/pure-ftpd.service
[Unit]
Description=Pure-FTPD daemon

[Service]
ExecStart=/usr/sbin/pure-config.pl /etc/pure-ftpd.conf
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process

[Install]
WantedBy=multi-user.target
