# octoprint-git/octoprint.service
[Unit]
Description=Octoprint 3d Printing Web Server
After=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStartPre=/usr/bin/mkdir -p /tmp/octoprint/ /srv/octoprint/timelapse/ /srv/octoprint/uploads/
ExecStart=/usr/bin/octoprint -c /etc/octoprint.conf --daemon start --pid /run/octoprint.pid
ExecStop=/usr/bin/octoprint --daemon stop --pid /run/octoprint.pid
StandardOutput=tty
StandardError=tty

[Install]
WantedBy=multi-user.target
