# fanout/fanout.service
[Unit]
Description=A simple pubsub message server"
After=network.target

[Service]
Type=forking
PIDFile=/run/fanout/fanout.pid
EnvironmentFile=/etc/conf.d/fanout
PermissionsStartOnly=true

ExecStartPre=/bin/mkdir -p /run/fanout
ExecStartPre=/bin/chown fanout:fanout /run/fanout
ExecStart=/usr/bin/fanout $FANOUT_OPTS
TimeoutSec=45s
Restart=always
User=fanout
Group=fanout

[Install]
WantedBy=multi-user.target
