# apache-activemq/activemq.service
[Unit]
Description=Popular and powerful open source messaging and Integration Patterns provider

[Service]
Type=forking
Environment=ACTIVEMQ_HOME=/opt/activemq
Environment=ACTIVEMQ_PIDFILE=/opt/activemq/data/activemq.pid
ExecStart=/opt/activemq/bin/activemq start
ExecStop=/opt/activemq/bin/activemq stop
ExecReload=/opt/activemq/bin/activemq restart
PIDFile=/opt/activemq/data/activemq.pid

[Install]
WantedBy=multi-user.target
