# tvheadend-git/tvheadend.service
[Unit]
Description=tvheadend
Wants=network.target



[Service]
EnvironmentFile=/etc/conf.d/tvheadend
PIDFile=/run/tvheadend.pid
ExecStart=/usr/bin/tvheadend $TVHEADEND_ARGS

[Install]
WantedBy=multi-user.target
