# tvheadend-git/tvheadend.service
[Unit]
Description=tvheadend
Wants=network.target

[Service]
User=hts
ExecStart=/usr/bin/tvheadend -C -u hts -g video

[Install]
WantedBy=multi-user.target
