# samsung-unified-driver-4.01.17/smfpd.service
[Unit]
Description=smfpd is a parallel port handling daemon.
After=network.target

[Service]
Type=forking
ExecStart=/usr/bin/smfpd

[Install]
WantedBy=multi-user.target