# shellinabox/shellinaboxd.service
[Unit]
Description=Serve a login-terminal over http on  port 4200.
Required=sshd.service
After=sshd.service

[Service]
User=root
Type=forking
EnvironmentFile=/etc/conf.d/shellinaboxd
ExecStart=/usr/bin/shellinaboxd  $SHELLINABOXD_ARGS
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=on-abort

[Install]
WantedBy=multi-user.target

