# battery_warn/battery_warn.service
[Unit]
Description=Blink Power LED when battery below threshold
After=acpid.service

[Service]
Type=simple
ExecStart=/usr/sbin/battery_warn
ExecReload=/bin/kill -9 $MAINPID
KillMode=process
Restart=always

[Install]
WantedBy=multi-user.target

