# radicale-git/radicale.service
[Unit]
Description=A Simple Calendar Server

[Service]
Type=forking
TimeoutStartSec=10
ExecStart=/usr/bin/radicale

[Install]
WantedBy=multi-user.target
