# i2p-bin/i2prouter.service
[Unit]
Description=I2P router
After=network.target

[Service]
ExecStart=/opt/i2p/i2prouter console
ExecReload=/opt/i2p/i2prouter restart
ExecStop=/opt/i2p/i2prouter stop

[Install]
WantedBy=multi-user.target
