# i2p-bin/i2prouter.service
[Unit]
Description=Invisible Internet Project
After=network.target

[Service]
Type=forking
PIDFile=/opt/i2p/.i2p/i2p.pid
ExecStart=/opt/i2p/i2prouter start
ExecReload=/opt/i2p/i2prouter restart
ExecStop=/opt/i2p/i2prouter graceful
TimeoutStopSec=5min

[Install]
WantedBy=multi-user.target
