# radicale/radicale.service
[Unit]
Description=RadiCAL Server
After=syslog.target network.target

[Service]
Type=forking
PIDFile=/run/radicale.pid
ExecStart=/usr/bin/radicale -d -p /run/radicale.pid

[Install]
WantedBy=multi-user.target
