# radicale/radicale.service
[Unit]
Description=RadiCAL Server

[Service]
ExecStart=/usr/bin/radicale
KillMode=process
