# radicale/radicale.service
[Unit]
Description=RadiCAL Server
After=syslog.target network.target

[Service]
Type=forking
ExecStart=/usr/bin/radicale
KillMode=process

[Install]
WantedBy=multi-user.target
