# mediabrowser-server-beta/mediabrowser-server.service
[Unit]
Description=Media Browser Server is a home media server built using other popular open source technologies.

[Service]
EnvironmentFile=/etc/conf.d/mediabrowser-server.conf
WorkingDirectory=/usr/lib/mediabrowser-server
ExecStart=/usr/bin/mediabrowser-server
User=mediabrowser
Group=mediabrowser
Restart=on-failure

[Unit]
After=network.target

[Install]
WantedBy=multi-user.target
