# tvheadend/tvheadend.service
[Unit]
Description=tvheadend
Wants=network.target,oscam.service
After=oscam.service

[Service]
User=hts
ExecStart=/usr/bin/tvheadend -C -u hts -g video

[Install]
WantedBy=multi-user.target
