# netkit-bsd-finger-ipv6/fingerd@.service
[Unit]
Description=Finger remote user information server

[Service]
ExecStart=/usr/sbin/in.fingerd
StandardInput=socket
