# wfrog-svn/wfrender.service
[Unit]
Description=wfrog renderer
After=network.target

[Service]
Type=simple
ExecStart=/usr/lib/wfrog/bin/wfrog --renderer -m
ExecReload=/usr/bin/kill -HUP $MAINPID
KillMode=process
Restart=always

[Install]
WantedBy=multi-user.target
