# prax-git/prax-system.service
[Unit]
Description=Prax System
# Note: target and service names are distribution-specific
# Requires=iptables.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/opt/prax/install/initd start
ExecStop=/opt/prax/install/initd stop
WorkingDirectory=/opt/prax

[Install]
WantedBy=prax-user.service
