# prax-git/prax-user.service
[Unit]
Description=Prax
#After=network.target
#Requires=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/opt/prax/bin/prax start
ExecStop=/opt/prax/bin/prax stop
WorkingDirectory=/opt/prax

[Install]
#WantedBy=default.target