# c-icap/c-icap.service
[Unit]
Description=ICAP server implementation

[Service]
EnvironmentFile=/etc/conf.d/c-icap
PIDFile=/var/run/c-icap/c-icap.pid
ExecStart=/usr/bin/c-icap $CICAP_ARGS

[Install]
WantedBy=multi-user.target
