# python2-leginon-svn/redux.service
[Unit]
Description=Use redux to serve images on myamiweb
Requires=network.target
After=syslog.target network.target systemd-user-sessions.service getty.target

[Service]
Type=simple
ExecStart=/usr/bin/reduxd
Restart=on-abort

[Install]
WantedBy=multi-user.target
