# policyd/policyd.service
[Unit]
Description=Policyd v2 daemon (ClueBringer)

[Service]
Type=forking
ExecStart=/usr/sbin/cbpolicyd --config /etc/policyd/cluebringer.conf
ExecStop=/bin/kill $MAINPID

[Install]
WantedBy=multi-user.target

