# kannel/bearerbox.service
[Unit]
Description=Kannel's Bearerbox


[Service]
SyslogIdentifier=bearerbox
PIDFile=/var/run/bearerbox.pid
ExecStart=/usr/bin/run_kannel_box --min-delay 5 --pidfile /var/run/bearerbox.pid /usr/bin/bearerbox /etc/kannel/kannel.conf

[Install]
WantedBy=multi-user.target
