# kannel/smsbox.service
[Unit]
Description=Kannel's SMSBox
Requires=bearerbox.service

[Service]
SyslogIdentifier=smsbox
PIDFile=/var/run/smsbox.pid
ExecStart=/usr/bin/run_kannel_box --min-delay 5 --pidfile /var/run/smsbox.pid /usr/bin/smsbox /etc/kannel/kannel.conf

[Install]
WantedBy=multi-user.target
