# kannel/wapbox.service
[Unit]
Description=Kannel's WAPBox
Requires=bearerbox.service

[Service]
SyslogIdentifier=wapbox
PIDFile=/var/run/wapbox.pid
ExecStart=/usr/bin/run_kannel_box --min-delay 5 --pidfile /var/run/wapbox.pid /usr/bin/wapbox /etc/kannel/kannel.conf

[Install]
WantedBy=multi-user.target
