# tpm-emulator/tpmd.service
[Unit]
Description=The TPM emulator daemon
After=sysinit.target

[Service]
User=tpm
ExecStart=/usr/bin/tpmd -f

[Install]
WantedBy=default.target
