# selenium-server-standalone/selenium-server.service
[Unit]
Description=Selenium Standalone Server

[Service]
EnvironmentFile=/etc/conf.d/selenium-server.conf
ExecStart=/usr/bin/selenium-server-java -jar /usr/share/selenium-server/selenium-server-standalone.jar ${SELENIUM_SERVER_ARGS}
KillMode=process
User=nobody

[Install]
WantedBy=multi-user.target
