# haproxy/haproxy.service
[Unit]
Description=HAProxy For TCP And HTTP Based Applications
After=network.target

[Service]
Type=forking
PIDFile=/run/haproxy.pid
ExecStartPre=/usr/sbin/haproxy -c -q -f /etc/haproxy/haproxy.cfg
ExecStart=/usr/sbin/haproxy -D -f /etc/haproxy/haproxy.cfg -p /run/haproxy.pid
ExecReload=-/bin/bash -c "exec /usr/sbin/haproxy -D -f /etc/haproxy/haproxy.cfg -p /run/haproxy.pid -sf $MAINPID"

[Install]
WantedBy=multi-user.target
