# blogspam/blogspam.service
[Unit]
Description=Blog-Spam XML-RPC service
After=network.target

[Service]
ExecStart=/usr/bin/vendor_perl/blogspam

[Install]
WantedBy=multi-user.target
