# mythtv-git-bas-t/mythbackend.service
[Unit]
Description=MythTV_backend
After=ffdecsawrapper.service mysqld.service network.target
Needs=ffdecsawrapper.service mysqld.service network.target

[Service]
Type=simple
Environment=HOME=/usr/share/mythtv
Nice=-4
User=root
#StandardOutput=null #(uncomment to disable mythbackend logs)
ExecStart=/usr/bin/mythbackend

[Install]
WantedBy=multi-user.target
