# riemann/riemann.service
[Unit]
Description=Riemann
After=network.target

[Service]
User=riemann
ExecStart=/usr/share/riemann/bin/riemann /etc/riemann/riemann.config

[Install]
WantedBy=multi-user.target
