# ovirt-vdsm-git/supervdsmd.service
[Unit]
Description="Auxiliary vdsm service for running helper functions as root"
Requires=libvirtd.service

[Service]
Type=simple
EnvironmentFile=-/etc/sysconfig/supervdsmd
ExecStart=/usr/share/vdsm/daemonAdapter "/usr/share/vdsm/supervdsmServer" --sockfile "/run/vdsm/svdsm.sock"
Restart=on-abort
