# prey-bash-client/prey-bash-client.timer
[Unit]
Description=Prey Bash Client timer
Documentation=file:/usr/share/prey/README

[Timer]
OnBootSec=1min
OnUnitActiveSec=20min
Unit=prey-bash-client.service

[Install]
WantedBy=multi-user.target
