# prey-bash-client/prey-bash-client.service
[Unit]
Description=Prey Tracker Bash Client
Documentation=file:/usr/share/prey/README

[Service]
Type=simple
ExecStart=/bin/sh -c '/usr/share/prey/prey.sh'
