# hadoop/hadoop-all.service
[Unit]
Description=Hadoop Datanode
Requires=hadoop-namenode.service hadoop-datanode.service hadoop-secondarynamenode.service hadoop-jobtracker.service hadoop-tasktracker.service

[Service]
User=hadoop
Group=hadoop

[Install]
WantedBy=multi-user.target
