# orientdb-systemd-unit/orientdb.service
[Unit]
Description=OrientDB NoSQL Graph-Document DBMS

[Service]
Type=simple
User=orient
SyslogIdentifier=orientdb

ExecStart=  /opt/orientdb/bin/server.sh \
	    1>>/opt/orientdb/log/orientdb.log \
	    2>>/opt/orientdb/log/orientdb.err
ExecStop=   /opt/orientdb/bin/shutdown.sh \
	    1>>/opt/orientdb/log/orientdb.log \
	    2>>/opt/orientdb/log/orientdb.err

[Install]
WantedBy=multi-user.target
