# trafmetux/trafmetux.cron
*/1  *  *  *  *  /usr/sbin/trafmetux update input output total
