# muhttpd/muhttpd.service
[Unit]
Description=Simple but complete web server
Documentation=man:muhttpd(8)
After=network.target

[Service]
Type=forking
ExecStart=/sbin/muhttpd

[Install]
WantedBy=multi-user.target
