# hhvm/hhvm@.service
[Unit]
Description=HHVM FastCGI application server
After=syslog.target network.target

[Service]
PIDFile=/run/hhvm/%I.pid
PrivateTmp=true
User=http
ExecStart=/usr/bin/hhvm --mode server -vPidFile=/run/hhvm/%I.pid --config=/etc/hhvm/%I.hdf --user http
Restart=on-failure

[Install]
WantedBy=multi-user.target