# conman/conman.service
[Unit]
Description=Conman Daemon
After=network.target

[Service]
Type=forking
ExecStart=/usr/sbin/conmand

[Install]
WantedBy=multi-user.target
