# icpld/icpld.service
[Unit]
Description=ICPLD
After=local-fs.target network.target

[Service]
Type=forking
ExecStart=/usr/bin/icpld
ExecStop=/usr/bin/icpld -quit

[Install]
WantedBy=multi-user.target
