# ers-node/ers-node.service
[Unit]
Description=ERS Node service
Requires=couchdb.service
After=couchdb.service

[Service]
User=ers-node
Type=simple
StandardOutput=journal
StandardError=journal
Restart=always
StartLimitInterval=10
StartLimitBurst=5
ExecStart=/usr/bin/python2 -m ers.daemon --pidfile=none --logtype=syslog

[Install]
WantedBy=multi-user.target

