# bdremote-ng/bdremote-ng.service
[Service]
EnvironmentFile=/etc/conf.d/bdremote-ng
Type=forking
ExecStart=/usr/sbin/bdremoteng -a ${REMOTE_ADDRESS}
