# kairosdb/kairosdb.service
[Unit]
Description=KairosDB is a fast distributed scalable time series database written on top of Cassandra.

[Service]
Type=forking
User=root
PIDFile=/var/run/kairosdb.pid
ExecStart=/opt/kairosdb/bin/kairosdb.sh run
ExecStop=/opt/kairosdb/bin/kairosdb.sh stop

[Install]
WantedBy=multi-user.target