# shadowsocks-nodejs-git/shadowsocks-nodejs.service
[Unit]
Description=A lightweight tunnel proxy.

[Service]
ExecStart=/usr/bin/sslocal /opt/shadowsocks-nodejs/config.json
StandardOutput=null

[Install]
WantedBy=multi-user.target
