# shadowsocks-nodejs-git/shadowsocks-nodejs.service
[Unit]
Description=A lightweight tunnel proxy.

[Service]
ExecStart=/usr/bin/node /opt/shadowsocks-nodejs/local.js
StandardOutput=null

[Install]
WantedBy=multi-user.target
