# python-nbviewer-git/nbviewer.service
[Unit]
Description=nbviewer webservice rendering ipynb to static HTML
 
[Service]
Type=simple
Restart=always
ExecStart=/usr/bin/python -m nbviewer --port=8080
 
[Install]
WantedBy=multi-user.target
