# greyhole/greyhole.service
[Unit]
Description=Greyhole
After=syslog.target network.target

[Service]
Type=simple
ExecStart=/usr/bin/greyhole -D
ExecReload=/usr/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
