# apt-cacher-ng/apt-cacher-ng.service
[Unit]
Description=A caching proxy specialized for package files
After=syslog.target network.target

[Service]
Type=forking
PIDFile=/run/apt-cacher-ng/pid
ExecStart=/usr/bin/apt-cacher-ng -c /etc/apt-cacher-ng/

[Install]
WantedBy=multi-user.target

