# sfptpd/sfptpd.service
[Unit]
Description=sfptpd Service
After=network.target
Requires=openonload.service
After=openonload.service

[Service]
Type=oneshot
RemainAfterExit=yes
PrivateTmp=true
EnvironmentFile=-/etc/default/sfptpd

ExecStart=/usr/bin/sfptpd -m ptp -f /etc/sfptpd.cfg
ExecStop=/usr/bin/killall -SIGINT /usr/bin/sfptpd

[Install]
WantedBy=multi-user.target
