# plymouth-git/kdm-unpatched-plymouth.service
[Unit]
Description=K Display Manager
After=systemd-user-sessions.service

Conflicts=plymouth-quit.service
After=plymouth-quit.service

[Service]
ExecStartPre=-/usr/bin/plymouth deactivate
ExecStartPre=-/usr/bin/plymouth quit --retain-splash
ExecStart=/usr/bin/kdm -nodaemon
StandardOutput=syslog

[Install]
Alias=display-manager.service
