# plymouth-git/kdm-unpatched-plymouth.service
[Unit]
Description=K Display Manager
After=systemd-user-sessions.service

Conflicts=getty@tty1.service plymouth-quit.service
After=getty@tty1.service plymouth-quit.service

[Service]
ExecStartPre=-/usr/bin/plymouth deactivate
ExecStartPre=-/usr/bin/plymouth quit --retain-splash
ExecStart=/usr/bin/kdm -nodaemon
StandardOutput=syslog

[Install]
WantedBy=graphical.target
