# iguana-ir/iguanair.service
[Unit]
Description=IguanaWorks infrared transceiver service
Documentation=https://github.com/terazini/abs-htpc/iguana-ir
Requires=local-fs.target

[Service]
EnvironmentFile=/etc/conf.d/iguana-ir
User=iguanair
Group=root
Type=forking
ExecStart=/usr/bin/igdaemon $DRIVERS --log-level=${LOGLEVEL} --send-timeout=${SENDTIMEOUT} --log-file=/tmp/igdaemon.log
ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/usr/bin/killall -w -s 2 igdaemon

[Install]
WantedBy=multi-user.target
