# chrontel-wetab/chrontel-wetab.service
[Unit]
Description=Chrontel HDMI Driver and Daemon
# Wait for the display manager if manual environment, e.g.
# Requires=gdm.service
# After=gdm.service

[Service]
# Manually force specific display, e.g. for lightdm you may use:
# Environment=DISPLAY=:0
# Environment=XAUTHORITY=/var/lib/lightdm/.Xauthority
ExecStart=/usr/sbin/ch7036_start

[Install]
WantedBy=graphical.target
