# agrarianskies2/agrarianskies2.service
[Unit]
Description=Minecraft Modpack Server Agrarian Skies 2

[Service]
ExecStart=/srv/agrarianskies2/agrarianskies2.sh start
ExecStop=/srv/agrarianskies2/agrarianskies2.sh stop
Type=forking

[Install]
WantedBy=multi-user.target
