# dm-cache-rootfs/dm-cache.service
[Unit]
Before=umount.target

[Install]
RequiredBy=local-fs.target

[Service]
Type=oneshot
ExecStart=/bin/true
ExecStop=/usr/lib/systemd/scripts/disable-dm-cache.sh
RemainAfterExit=yes
