# sickrage/sickrage-system.service
[Unit]
Description=SickRage Daemon
After=network.target

[Service]
User=sickrage
Group=sickrage
ExecStart=/usr/bin/env python2 /opt/sickrage/SickBeard.py --quiet --config /opt/sickrage/config.ini --datadir /opt/sickrage

[Install]
WantedBy=multi-user.target
