# systemd-rtirq/rtirq.service
[Unit] 
Description=Realtime IRQ thread system tuning 

[Service] 
Type=oneshot 
ExecStart=/usr/lib/systemd/scripts/rtirq start
ExecStop=/usr/lib/systemd/scripts/rtirq stop 
RemainAfterExit=true 

[Install] 
WantedBy=multi-user.target
