# happy-hacking-linux/happy-hacking-linux.service
[Unit]
Description=Happy Hacking Linux
# After=systemd-user-sessions.service

[Service]
ExecStart=/home/anekos/project/happy-hacking-linux/source/happy-hacking-linux

[Install]
WantedBy=multi-user.target
