# timer_entropyd/timer_entropyd.service
[Unit]
Description=Gathering entropy from system timers

[Service]
Type=forking
ExecStart=/usr/bin/timer_entropyd

[Install]
WantedBy=multi-user.target
