# h2oai/h2oai.service
[Unit]
Description=H2O Open source big data analysis tool
After=network.target

[Service]
EnvironmentFile=/etc/conf.d/h2oai
ExecStart=-/usr/bin/java -jar /usr/share/java/h2o/h2o.jar $ARGS
Type=simple
Restart=on-failure
RestartSec=5
StartLimitInterval=60s
StartLimitBurst=3
StandardOutput=null

[Install]
WantedBy=multi-user.target
