# sogo/sogo.service
[Unit]
Description=Sogo Daemon
After=network.target slapd.service memcached.service postfix.service dovecot.service postgresql.service

[Service]
Type=forking
User=sogo
ExecStart=/usr/sbin/sogod -WOPidFile /var/run/sogo/sogo.pid
PIDFile=/var/run/sogo/sogo.pid

[Install]
WantedBy=multi-user.target
