# freedns-daemon/freedns-daemon.service
[Unit]
Description=Updating freedns record
After=syslog.target network.target

[Service]
Type=simple
ExecStart=/usr/bin/freedns-daemon

[Install]
WantedBy=multi-user.target
