# tor-git/tor.service
[Unit]
Description=Anonymizing overlay network
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/tor --runasdaemon 0 --quiet
ExecReload=/bin/kill -HUP ${MAINPID}
ExecStop=/bin/kill -INT ${MAINPID}
TimeoutSec=30
Restart=on-failure
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target
