# pcsc-cyberjack/99-cyberjack.rules
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="0c4b", \
  GROUP="pcscd", MODE="660", \
  RUN+="/bin/mkdir /var/run/pcscd", \
  RUN+="/bin/chgrp pcscd /var/run/pcscd", \
  RUN+="/bin/chmod 0775 /var/run/pcscd"