# apparmor/apparmor.service
[Unit]
Description=AppArmor profiles

[Service]
Type=oneshot
ExecStart=/usr/bin/apparmor_load.sh
ExecStop=/usr/bin/apparmor_unload.sh
RemainAfterExit=yes

[Install]
WantedBy=basic.target
