# p4d/p4d.service
[Unit]
Description=Perforce Server
Requires=network.target local-fs.target

[Service]
Type=forking
ExecStart=/usr/bin/p4d -d -r /home/p4/
ExecStop=/usr/bin/killall p4d
WorkingDirectory=/home/p4/

[Install]
WantedBy=multi-user.target
