# wildfly-devel/wildfly-devel.service
[Unit]
Description=Wildfly Application Server

[Service]
Type=simple
User=wildfly
PIDFile=/run/wildfly-devel/wildfly-devel.pid
ExecStart=/usr/share/wildfly-devel/bin/systemd-wrapper.sh
ExecReload=/usr/share/wildfly-devel/bin/jboss-cli.sh --connect command=:reload
ExecStop=/usr/share/wildfly-devel/bin/jboss-cli.sh --connect command=:shutdown

[Install]
WantedBy=multi-user.target
