# rabbitmq-systemd-unit/rabbitmq.service
[Unit]
Description=RabbitMQ Messaging Server

[Service]
Type=forking
User=rabbitmq
SyslogIdentifier=rabbitmq
EnvironmentFile=/etc/rabbitmq/rabbitmq-env.conf
PIDFile=$MNESIA_HOME/rabbit.pid

#ExecStartPre=
ExecStart=  /usr/sbin/rabbitmq-server -detached
ExecStop=   /usr/sbin/rabbitmqctl stop

# OOMScoreAdjust=-200

[Install]
WantedBy=multi-user.target
