# rabbitmq-systemd-unit/rabbitmq.service
[Unit]
Description=RabbitMQ Messaging Server

[Service]
Type=forking
SyslogIdentifier=rabbitmq
EnvironmentFile=/etc/rabbitmq/rabbitmq-env.conf
PIDFile=$MNESIA_HOME/rabbit.pid

#ExecStartPre=
ExecStart=  /bin/su - rabbitmq -m -c "/usr/sbin/rabbitmq-server -detached"
ExecReload= /usr/sbin/rabbitmqctl stop && /bin/su - rabbitmq -m -c "/usr/sbin/rabbitmq-server -detached"
ExecStop=   /usr/sbin/rabbitmqctl stop

# OOMScoreAdjust=-200

[Install]
WantedBy=multi-user.target
