# distcc-raspi/distccd-raspi.service
[Unit]
Description=A distributed C/C++ compiler for RaspberryPi

[Service]
EnvironmentFile=/etc/conf.d/distccd-raspi.conf
ExecStart=/usr/bin/distccd-raspi --daemon --no-detach $DISTCC_ARGS
PrivateTmp=true

[Install]
WantedBy=multi-user.target
