# ramlog/ramlog.service
[Unit]
Description=Log Storage RamDisk

[Service]
ExecStart=/etc/ramlog/ramlog start > /dev/null
ExecStop=/etc/ramlog/ramlog stop > /dev/null
Type=oneshot
RemainAfterExit=yes
After=local-fs.target

[Install]
WantedBy=multi-user.target
