# karmacoin-daemon/karmacoind@.service
[Unit]
Description=DogeCoin Daemon
After=network.target 

[Service]
User=%i
Group=%i
Type=simple
ExecStart=/usr/bin/karmacoind

[Install]
WantedBy=multi-user.target
