# wondershaper-git/wondershaper.service
[Unit]
Description=Bandwidth shaper/Network rate limiter
After=network.target
Wants=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/wondershaper -a eth0 -d 2048 -u 512
ExecStop=/usr/bin/wondershaper -c -a eth0

[Install]
WantedBy=multi-user.target