# tor-control-port-proxy/tor-control-port-proxy.service
[Unit]
Description=Tor control port proxy
After=tor.service

[Service]
User=tor
Type=simple
ExecStart=/usr/bin/tor-control-port-proxy

[Install]
WantedBy=multi-user.target
