# logrotate-hourly/logrotate.timer
[Unit]
Description=Hourly rotation of log files

[Timer]
OnCalendar=hourly
AccuracySec=1h
Persistent=true
