# cloudprint-daemon-git/cloudprint-daemon-git.service
[Unit]
Description=Share your CUPS printers with google's cloud print.
Requires=network.target cups.service
After=network.target cups.service NetworkManager-wait-online.service

[Service]
Type=simple
ExecStart=/usr/bin/cloudprint
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
