# collectw/collectw.service
[Unit]
Description=Web-based UI for collectd.

[Service]
Type=simple
User=nobody
ExecStart=/usr/bin/collectw

[Install]
WantedBy=multi-user.target
