# arch-macbookair6.2-config/systemd-warnings.install
post_install() {


cat << EOF
It's recommended to execute the following commands to enable the systemd services:
systemctl enable gpe66-fix.service
systemctl enable powertop.service
EOF

}