# tinc-pre-systemd/tincd@.service
[Unit]
Description=tinc vpn on network %i
After=network.target

[Service]
Type=forking
ExecStart=/usr/sbin/tincd --user=nobody --chroot --mlock --net=%i

[Install]
WantedBy=multi-user.target
