# webogram-git/webogram.service
[Unit]
Description=Telegram unofficial web app

[Service]
WorkingDirectory=/opt/webogram
ExecStart=/usr/bin/node server.js 3772
Restart=always

[Install]
WantedBy=multi-user.target
