# riak-bin/riak.service
[Unit]
Description=Distributed key/value store from Basho Technologies
After=network.target

[Service]
User=riak
Type=forking
ExecStart=/usr/sbin/riak start
ExecStop=/usr/sbin/riak stop
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target
