# iodine-git/iodined.service
[Unit]
After=network.target
Requires=network.target

[Service]
EnvironmentFile=/etc/conf.d/iodined
ExecStart=/usr/sbin/iodined $IODINED_ARGS

[Install]
WantedBy=multi-user.target

