# iodine-git/iodine.service
[Unit]
After=network.target
Requires=network.target

[Service]
EnvironmentFile=/etc/conf.d/iodine
ExecStart=/usr/sbin/iodine $IODINE_ARGS

[Install]
WantedBy=multi-user.target

