# hiveary-agent/hiveary-agent.service
[Unit]
Description=An intelligent server monitoring solution that automatically finds anomalies.
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/hiveary-agent
TimeoutStopSec=20
KillSignal=SIGINT
Restart=on-abort

[Install]
WantedBy=multi-user.target
