# driveclient/driveclient.service
[Unit]
Description=Rackspace backup agent
After=network.target

[Service]
ExecStart=/usr/bin/driveclient

[Install]
WantedBy=multi-user.target
