# driveclient/driveclient.service
[Unit]
Description=Rackspace backup agent
After=network.target

[Service]
Type=simple
Environment= ARGS=-d
EnvironmentFile=-/etc/conf.d/driveclient
ExecStart=/usr/local/bin/driveclient $ARGS

[Install]
WantedBy=multi-user.target
