# greylistd/greylistd.service
[Unit]
Description=Greylist Daemon

[Service]
User=greylist
ExecStart=/usr/sbin/greylistd
ExecReload=/usr/bin/greylist reload

[Install]
WantedBy=multi-user.target
