# cube-git/cube-evaluator.service
[Unit]
Description=Cube Evaluator
Requires=mongodb.service
After=mongodb.service

[Service]
User=cube
Group=cube
WorkingDirectory=/usr/share/cube
ExecStart=/usr/bin/node bin/evaluator.js 2>&1 >> /var/log/cube/evaluator.log

[Install]
WantedBy=multi-user.target
