# systemd-ckboot/ckboot.service
[Unit]
Description=Checks boot files and MBR for changes
After=local-fs.target

[Service]
Type=simple
ExecStart=/usr/bin/ckboot-check
