# shadow-systemd/verify-shadow.service
[Unit]
Description=Verify integrity of password and group files

[Service]
Type=oneshot
ExecStart=/usr/sbin/pwck -r
ExecStart=/usr/sbin/grpck -r

[Install]
WantedBy=cron-daily.target
