# vnstatui/vnstatui.service
[Unit]
Description=VNStat Web UI
Requires=vnstat.service
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/bin/uwsgi --ini /usr/lib/python3.3/site-packages/vnstatui/vnstatui.ini
WorkingDirectory=/usr/lib/python3.3/site-packages/vnstatui

[Install]
WantedBy=multi-user.target
