# vnstatui/vnstatui.service
[Unit]
Description=VNStat Web UI
Requisite=vnstat.service
Requires=vnstat.service
After=vnstat.service
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/bin/vnstatui
WorkingDirectory=/usr/lib/python3.3/site-packages/vnstatui

[Install]
WantedBy=multi-user.target
