# vnstatui/vnstatui.service
[Unit]
Description=VNStat Web UI
Requires=vnstat.service
Wants=network.target

[Service]
EnvironmentFile=/etc/conf.d/vnstatui.uwsgi
ExecStart=/usr/bin/uwsgi --pidfile /run/uwsgi-vnstatui.pid
ExecStop=/usr/bin/uwsgi --stop /run/uwsgi-vnstatui.pid
WorkingDirectory=/usr/lib/python3.4/site-packages/vnstatui

[Install]
WantedBy=multi-user.target
