# ppcoind/ppcoind.service
[Unit]
Description=PPCoin Daemon
After=network.target

[Install]
WantedBy=multi-user.target

[Service]
User=ppcoind

Type=forking
TimeoutStartSec=30
PIDFile=/run/ppcoind/ppcoind.pid
ExecStart=/usr/bin/ppcoind -daemon \
                           -conf=/etc/ppcoind/ppcoind.conf \
                           -pid=/run/ppcoind/ppcoind.pid \
                           -datadir=/var/lib/ppcoind
