# zfoned/zfoned.service
[Unit]
Description=Zfone Daemon
After=network.target

[Service]
#User=
ExecStart=/usr/local/bin/zfoned
ExecStop=/bin/kill -INT $MAINPID

[Install]
WantedBy=multi-user.target
