# glassfish-v3/glassfish.service
[Unit]
Description = GlassFish Server
After = syslog.target network.target

[Service]
ExecStart = /opt/glassfish3/glassfish/bin/asadmin start-domain
ExecStop = /opt/glassfish3/glassfish/bin/asadmin stop-domain
ExecReload = /opt/glassfish3/glassfish/bin/asadmin restart-domain
Type = forking

[Install]
WantedBy = multi-user.target
