# python2-pifacenetctl/pifacenetctl.service
[Unit]
Description=PiFace headless monitor/control for wireless networks
After=network.target

[Service]
Type=forking
PIDFile=/var/run/pifacenetctl.pid
ExecStart=/usr/bin/python2 -m pifacenetctl start
ExecStop=/usr/bin/python2 -m pifacenetctl stop

[Install]
WantedBy=multi-user.target

