# btlive/btlive.service
[Unit]
Description=BTLive Daemon
After=network.target

[Service]
#User=
ExecStart=/usr/bin/btlive
ExecStop=/bin/kill -INT $MAINPID

[Install]
WantedBy=multi-user.target
