# feedthebeast-server/feedthebeastd.service
[Unit]
Description=Feed the Beast server daemon
After=network.target

[Service]
ExecStart=/usr/bin/feedthebeastd
Type=oneshot
RemainAfterExit=yes
ExecStop=/usr/bin/screen -S feedthebeast -X stuff 'say Server Shutting Down in 5 Seconds\n' ; /usr/bin/sleep 5 ; /usr/bin/screen -S feedthebeast -X stuff 'stop\n' ; /usr/bin/expect -c "exec screen -x feedthebeast ; wait ; exit"

[Install]
WantedBy=multi-user.target
