# btrbk/btrbk.service
[Unit]
Description = BTRBK Backup

[Service]
Type=simple
ExecStart=/usr/bin/btrbk run
