# etherpad-lite/etherpad-lite.service
[Unit]
Description=Etherpad-Lite server
After=syslog.target
After=network.target

[Service]
Type=simple
User=etherpad
Group=etherpad
PermissionsStartOnly=true
ExecStartPre=-/usr/bin/mkdir /var/run/etherpad-lite
ExecStartPre=/usr/bin/chown -R etherpad:etherpad /var/run/etherpad-lite
ExecStart=/usr/share/webapps/etherpad-lite/bin/run.sh

[Install]
WantedBy=multi-user.target
