# etherpad-lite/etherpad-lite.service
[Unit]
Description=Etherpad-Lite server
After=syslog.target
After=network.target

[Service]
Type=simple
User=etherpad
Group=etherpad
ExecStart=/usr/share/webapps/etherpad-lite/bin/run.sh

[Install]
WantedBy=multi-user.target
