# ndppd/ndppd.service
[Unit]
Description="ndppd NDP Proxy Daemon"
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/bin/ndppd
StandardInput=null
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
