# ptpd2/ptpd2.service
[Unit]
Description=Precision Time Protocol daemon
After=network.target

[Service]
Type=forking
ExecStart=/usr/sbin/ptpd2 -c /etc/ptpd2.conf
PIDFile=/var/run/ptpd2.lock

[Install]
WantedBy=multi-user.target
