# haste-server-git/haste-server.service
[Unit]
Description=Haste Server
Requires=redis.service

[Service]
WorkingDirectory=/usr/lib/node_modules/haste
ExecStart=/usr/bin/node /usr/bin/haste-server

[Install]
WantedBy=multi-user.target
