# bogom/bogom.service
[Unit]
Description=A simple sendmail milter to interface bogofilter.
After=network.target

[Service]
User=bogom
Group=mail
UMask=0002
Type=forking
PIDFile=/run/bogom/bogom.pid
ExecStart=/usr/bin/bogom -c /etc/bogom.conf

[Install]
WantedBy=multi-user.target
