# sssd/sssd.service
[Unit]
Description=System Security Services Daemon
Before=rpc-idmapd.service

[Service]
ExecStart=/usr/sbin/sssd -D -f
Type=forking
PIDFile=/run/sssd.pid

[Install]
WantedBy=multi-user.target
