# connman-dispatcher/connman-dispatcher.service
[Unit]
Description=Network dispatcher service
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/connman-dispatcher -p /etc/connman-dispatcher/

[Install]
WantedBy=multi-user.target
