# procwd/procwd.service
[Unit]
Description=A utility for monitoring processes, restarting them if they die

[Service]
Type=forking
ExecStart=/usr/bin/procwd

[Install]
WantedBy=multi-user.target

