# freefall/freefall@.service
[Unit]
Description=Disk protection for HP/DELL machines.

[Service]
Type=forking
# freefall takes disk device as argument.
ExecStart=/usr/bin/freefall "%I"

[Install]
# The instance name is escaped.
DefaultInstance=-dev-sda
WantedBy=multi-user.target

