# cndrvcups-capt/ccpd.service
[Unit]
Description=Canon CAPT Printing Daemon

[Service]
ExecStart=/usr/sbin/ccpd
Type=forking

[Install]
WantedBy=multi-user.target
