# abyssws/abyssws.service
[Unit]
Description=Abyss Web Server X1

[Service]
EnvironmentFile=/etc/conf.d/abyssws.conf.d
ExecStart=/opt/abyssws/abyssws -d $OPTIONS
Type=forking

[Install]
WantedBy=multi-user.target