# headphones-git/headphones.service
[Unit]
Description=Headphones Daemon

[Service]
EnvironmentFile=/etc/conf.d/headphones
ExecStart=/usr/bin/headphones -q -d --config ${HP_CONF} --datadir ${HP_DIR} --port ${HP_PORT}
GuessMainPID=no
Type=forking
User=headphones
Group=headphones

[Install]
WantedBy=multi-user.target
