# systemd-timer-pkgcacheclean/pkgcacheclean3.service
[Unit]
Description=Remove all but the last two versions of pkgs in the pacman cache.

[Service]
Type=oneshot
ExecStart=/usr/bin/pkgcacheclean 3

