# systemd-timer-pkgcacheclean/pkgcacheclean.timer
[Unit]
Description=Remove all but the last two versions of pkgs in the pacman cache, every month.

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target

