# flappycoin-daemon/flappycoind@.service
[Unit]
Description=FlappyCoin Daemon
After=network.target 

[Service]
User=%i
Group=users
Type=simple
ExecStart=/usr/bin/flappcoind

[Install]
WantedBy=multi-user.target
