# systemd-vgaswitcheroo-units/vgaswitcheroo.service
[Unit]
Description=Deactivate the discrete graphic

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/vgaswitcheroo_start.sh
ExecStop=/usr/bin/vgaswitcheroo_stop.sh

[Install]
WantedBy=multi-user.target
