# systemd-vgaswitcheroo-units/vgaswitcheroo.service
[Unit]
Description=Deactivate the discrete graphic

[Service]
ExecStart=/usr/bin/vgaswitcheroo.sh

[Install]
WantedBy=multi-user.target
