# remotefs/rfsd.service
[Unit]
Description=Remote Filesystem daemon
After=syslog.target network.target

[Service]
Type=forking
PIDFile=/var/run/rfsd.pid
EnvironmentFile=/etc/default/rfsd
ExecStart=/usr/bin/rfsd -a $RFSD_LISTEN -p $RFSD_PORT -u $RFSD_USER -r $RFSD_PIDFILE $RFSD_QUIET
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
