# athcool/athcool.service
[Unit]
Description=Athcool
Wants=local-fs.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/sbin/athcool on
ExecStop=/usr/sbin/athcool off

[Install]
WantedBy=multi-user.target
