# foldingathome-smp/foldingathome-smp.service
[Unit]
Description=Folding@Home distributed computing client
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/fah-smp
ExecStart=/opt/fah-smp/fah6 -smp -verbosity 9 -forceasm 

[Install]
WantedBy=multi-user.target
