# xbmc-eventclients-wiiremote/xbmc-wiiremote.service
[Unit]
Description=XBMC Media Center Wiimote (Wii Remote) Eventclient
Documentation=man:xbmc-wiiremote(1) http://wiki.xbmc.org/index.php?title=EventServer#WiiRemote

[Service]
ExecStart=/usr/bin/xbmc-wiiremote

[Install]
WantedBy=multi-user.target
