# petrified/petrified.crontab
## Various examples of how to run petrified via cron -- note that not all
## cron implementations support /etc/cron.d/ (i.e. fcron does not)

## run system level every 15 minutes
## /etc/cron.d/petrified
*/15 * * * * root /usr/bin/petrified

## run system level every hour
## /etc/cron.hourly/petrified
#!/bin/bash
/usr/bin/petrified

## run user level every 30 minutes
## (f)crontab -e (-u <username>)
*/30 * * * * /usr/bin/petrified
