# irtrans-irserver/irtrans-irserver.service
[Unit]
Description=IRTrans Server
BindsTo=dev-$IRDEVICE.device

[Service]
Type=forking
PIDFile=/run/irtrans-irserver.pid
EnvironmentFile=/etc/conf.d/irtrans-irserver
ExecStart=/usr/bin/irserver $IRSERVER_ARGS -daemon -logfile $LOGFILE -loglevel $LOGLEVEL -pidfile /run/irtrans-irserver.pid $IRDEVICE

[Install]
WantedBy=multi-user.target
