# novacom-systemd/novacom.service
[Unit]
Description=Novacom service for webOS development connectivity

[Service]
Type=forking
EnvironmentFile=/etc/conf.d/novacomd
GuessMainPID=yes
ExecStart=/usr/bin/novacomd $NOVACOMD_OPTIONS

[Install]
WantedBy=multi-user.target
