# burp-backup-dev/burp-client.service
[Unit]
Description=Burp Backup Script

[Service]
Type=simple
ExecStart=/usr/sbin/burp -a t
