# wlan-powersave/70-wifi-powersave.rules
# Automatically turn on power save for wireless devices
ACTION=="add", SUBSYSTEM=="net", KERNEL=="wlan*", RUN+="/usr/bin/iw dev %k set power_save on"
