# ircd-hybrid/ircd-hybrid.service
[Unit]
Description=IRC daemon
After=network.target remote-fs.target nss-lookup.target

[Service]
Type=forking
User=ircd
PIDFile=/var/run/ircd/ircd-hybrid.pid
ExecStart=/usr/sbin/ircd -logfile /var/log/ircd/ircd.log -configfile /etc/ircd-hybrid/ircd.conf

[Install]
WantedBy=multi-user.target
