# rdiffweb/rdiffweb.service
[Unit]
Description=Web interface for rdiff-backup
After=network.target

[Service]
Type=forking
ExecStart=/usr/bin/rdiff-web --pid-file=/run/rdiffweb.pid -d --background 
ExecStartPre=/bin/rm -f /run/rdiffweb.pid
PIDFile=/run/rdiffweb.pid
