# triggerhappy-git/triggerhappy.service
[Unit]
Description=triggerhappy hotkey daemon

[Service]
Type=simple
ExecStart=/usr/bin/thd --triggers /etc/triggerhappy.conf --socket /var/run/thd.socket /dev/input/event*

[Install]
WantedBy=multi-user.target
