# systemd-redis-units/redis.service
[Unit]
Description=Redis
After=network.target

[Service]
Type=forking
ExecStart=/usr/bin/redis-server /etc/redis.conf
Restart=always

[Install]
WantedBy=multi-user.target
