# ccnx/ndnd.service
[Unit]
Description=Content-Centric Networking forwarder and router
After=network.target

[Service]
Type=forking
ExecStart=/usr/bin/ccndstart
ExecReload=/usr/bin/ccndstop && /usr/bin/ccndstart
ExecStop=/usr/bin/ccndstop

[Install]
WantedBy=multi-user.target
