# fmberry-rpi-git/fmberry.service
[Unit]
Description = Start the FMBerry daemon

[Service]
ExecStart=/usr/bin/fmberryd
Restart=always

[Install]
WantedBy = multi-user.target