# oftpd/oftpd.service
[Unit]
Description=oftpd ftp server
After=network.target

[Service]
Type=simple
EnvironmentFile=/etc/conf.d/oftpd
ExecStart=/usr/bin/oftpd -N $OFTPD_ARGS "$OFTPD_USER" "$OFTPD_PATH"
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
