# razercfg/razerd.service
[Unit]
Description=Razer device state daemon

[Service]
ExecStart=/usr/sbin/razerd

[Install]
WantedBy=graphical.target
