# open-sasc-ng/sasc-ng.service
[Unit]
Description=A SoftCAM creating virtual DVB interface
After=oscam.service 
Needs=oscam.service

[Service]
Type=simple
User=root
Nice=-6
Restart=Restart-always
EnvironmentFile=/etc/conf.d/sasc-ng
ExecStartPre=/usr/bin/modprobe dvbloopback $DVBLOOPBACK_ARGS
TimeoutSec=2
ExecStart=/usr/bin/sasc-ng $SASCNG_ARGS --cam-dir $CAMDIR --log $LOGFILE

[Install]
WantedBy=multi-user.target
