# open-sasc-ng/sasc-ng.service
[Unit]
Description=Sasc-ng
After=oscam.service

[Service]
Type=simple
EnvironmentFile=/etc/conf.d/sasc-ng
ExecStartPre=/sbin/modprobe dvbloopback $DVBLOOPBACK_ARGS
TimeoutSec=2
ExecStart=/usr/sbin/sasc-ng $SASCNG_ARGS --cam-dir $CAMDIR --log $LOGFILE
TimeoutSec=2

[Install]
WantedBy=multi-user.target
