# open-sasc-ng/sasc-ng.service
[Unit]
Description=Sasc-ng
After=network.target oscam.service

[Service]
Type=forking
EnvironmentFile=/etc/conf.d/sasc-ng.conf
ExecStartPre=/sbin/modprobe dvbloopback $DVBLOOPBACK_ARGS
TimeoutSec=2
ExecStart=/usr/sbin/sasc-ng -D $SASCNG_ARGS --cam-dir=$CAMDIR -l $LOGDIR/sasc-ng.log

[Install]
WantedBy=multi-user.target
