# open-sasc-ng/sasc-ng.service
[Unit]
Description=Sasc-ng
After=oscam.service
Needs=oscam.service

[Service]
Type=simple
Restart=Restart-always
EnvironmentFile=/etc/conf.d/sasc-ng
ExecStartPre=/sbin/modprobe dvbloopback $DVBLOOPBACK_ARGS
TimeoutSec=2
ExecStart=/usr/bin/sasc-ng $SASCNG_ARGS --cam-dir $CAMDIR --log $LOGFILE

[Install]
WantedBy=multi-user.target
