# open-sasc-ng/sasc-ng.service
[Unit]
Description=Sasc-ng
After=oscam.service

[Service]
Type=forking
EnvironmentFile=/etc/conf.d/sasc-ng
ExecStartPre=/sbin/modprobe dvbloopback $DVBLOOPBACK_ARGS
TimeoutSec=2
ExecStart=/usr/sbin/sasc-ng -D $SASCNG_ARGS --cam-dir=$CAMDIR
TimeoutSec=4

[Install]
WantedBy=multi-user.target
