# buildhck-git/buildhck.service
[Unit]
Description=buildhck service unit
Requires=uwsgi@buildhck.service

[Service]
Type=oneshot
ExecReload=/usr/bin/systemctl reload uwsgi@buildhck.service
ExecStop=/usr/bin/systemctl stop uwsgi@buildhck.service
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
