# aiccu/aiccu.service
[Unit]
Description=AICCU (Automatic IPv6 Connectivity Configuration Utility)
Requires=network.target
After=syslog.target network.target time-sync.target 

[Service]
Type=forking
PIDFile=/var/run/aiccu.pid
ExecStart=/usr/sbin/aiccu start $OPTIONS
ExecStop=/usr/sbin/aiccu stop

[Install]
WantedBy=multi-user.target
