# teamspeak3-server/teamspeak3-server.service
[Unit]
Description=TeamSpeak3 Server
After=network.target

[Service]
Type=simple
ExecStart=/bin/su -s /bin/bash -l -c "/usr/bin/teamspeak3-server logpath=/var/log/teamspeak3-server/ dbsqlpath=/usr/share/teamspeak3-server/sql/" teamspeak

[Install]
WantedBy=multi-user.target
