# multigl-libgl/multigl.service
[Unit]
Description=Handles kernel commandline setting multigl.
Before=graphical.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/switchgl --init

[Install]
WantedBy=graphical.target