# prey-tracker/prey-tracker.service
[Unit]
Description=Prey Tracker Service
Documentation=/usr/share/prey/README

[Service]
Type=simple
ExecStart=/usr/share/prey/prey.sh > /var/log/prey.log
