# plymouth-release/gdm-plymouth.service
[Unit]
Description=GNOME Display Manager
After=systemd-user-sessions.service

Conflicts=plymouth-quit.service
After=plymouth-quit.service

[Service]
ExecStart=/usr/sbin/gdm -nodaemon
StandardOutput=syslog

[Install]
Alias=display-manager.service
