# plymouth-release/lxdm-plymouth.service
[Unit]
Description=LXDE Display Manager
After=systemd-user-sessions.service

Conflicts=plymouth-quit.service
After=plymouth-quit.service

[Service]
ExecStart=/usr/sbin/lxdm >/dev/null 2>&1 </dev/null
StandardOutput=syslog

[Install]
Alias=display-manager.service
