# alljoyn/alljoyn-daemon.service
[Unit]
Description=Alljoyn Router Daemon
After=network.target

[Service]
ExecStart=/usr/bin/alljoyn-daemon --verbosity=3 --nofork

[Install]
WantedBy=multi-user.target
