# alljoyn/alljoyn-daemon.socket
[Unit]
Description=Alljoyn socket for unix connection
PartOf=alljoyn-daemon.service

[Socket]
ListenStream=@alljoyn
ListenStream=0.0.0.0:9955
ListenDatagram=0.0.0.0:9955
SocketMode=0666

[Install]
WantedBy=sockets.target
