# shxd/shxd.service
[Unit]
Description=Synapse Hotline X Daemon
After=network.target

[Service]
ExecStart=/srv/shxd/hxd
Type=oneshot
RemainAfterExit=yes
ExecStop=kill hxd

[Install]
WantedBy=multi-user.target
 