# wuala-daemon/wuala.service
[Unit]
Description=Wuala as service
After=network.target

[Service]
EnvironmentFile=/etc/conf.d/wualad
ExecStart=/bin/su --command "wuala_login" $LinuxLoginName
ExecStop=/bin/su --command "~/wuala/wuala shutdown" $LinuxLoginName

[Install]
WantedBy=multi-user.target
