# pandora-libraries-git/pndnotifyd.service
[Unit]
Description=Manages the auto-discovery mechanism and integrating PXML.xml dir-apps and .pnd files into the FreeDesktop standard system (.desktop files).
After=local-fs.target

[Service]
ExecStart=/usr/bin/pndnotifyd -d
Type=forking

[Install]
WantedBy=multi-user.target
