# liquidsoap-git/liquidsoap-git.service
[Unit]
Description=Liquidsoap service for channel %I
Documentation=man:liquidsoap(1)
After=network.target
ConditionPathExists=/etc/liquidsoap/%I.liq

[Service]
Type=forking
PIDFile=/run/liquidsoap/%I.pid
ExecStart=/usr/bin/liquidsoap --daemon /etc/liquidsoap/%I.liq
Restart=on-failure
User=liquidsoap
Group=liquidsoap

[Install]
WantedBy=multi-user.target

