# update-pacman-mirrorlist/update-pacman-mirrorlist.timer
[Unit]
Description=Updates the pacman mirrorlist from the Pacman Mirrorlist Generator on a weekly basis.

[Timer]
OnCalendar=weekly
AccuracySec=1d
Persistent=true
