# update-pacman-mirrorlist/update-pacman-mirrorlist.service
[Unit]
Description=Updates the pacman mirrorlist from the Pacman Mirrorlist Generator.
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/update-pacman-mirrorlist '/etc/pacman.d/mirrorlist' 'https://www.archlinux.org/mirrorlist/?country=all&use_mirror_status=on'
