# update-pacman-mirrorlist/update-pacman-mirrorlist.timer
[Unit]
Description=Updates the pacman mirrorlist on a weekly basis.

[Timer]
OnCalendar=weekly
AccuracySec=1d
Persistent=true
