# ntopng-git/ntopng.service
[Unit]
Description=ntopng network traffic monitor

[Service]
ExecStart=/usr/bin/ntopng

[Install]
WantedBy=multi-user.target
