# daemontools/svscan.service
[Unit]
Description=svscan service of daemontools 
After=syslog.target

[Service]
ExecStart=/usr/sbin/svscanboot
Restart=always

[Install]
WantedBy=multi-user.target
