# opensand-collector/opensand-collector@.service
[Unit]
Description="Statistics and events collector for OpenSAND"
After=syslog.target network.target

[Service]
User=opensand
ExecStart=sand_collector -i %I

[Install]
WantedBy=multi-user.target

