# rspamd/rspamd.service
[Unit]
Description=Fast, free and open-source spam filtering system.
After=network.target

[Service]
ExecStartPre=/usr/bin/rspamd -t -c /etc/rspamd.xml -u rspamd -g rspamd
ExecStart=/usr/bin/rspamd -f -c /etc/rspamd.xml -u rspamd -g rspamd
ExecReload=/usr/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
