# ghost/ghost.service
[Unit]
Description=Ghost Blogging Platform
After=network.target

[Service]
ExecStart=/usr/bin/npm start --production /srv/ghost

[Install]
WantedBy=multi-user.target
