# cras-daisy-spring/cras.socket
[Unit]
Description=Chromium OS audio server socket

[Socket]
ListenStream=/var/run/cras/.cras_socket

[Install]
WantedBy=cras.service
