# cras-daisy-spring/cras.service
[Unit]
Description=ChromeOS audio server
Requires=dbus.service cras.socket cras-state.service
After=dbus.service cras.socket
Before=cras-state.service

[Service]
Type=dbus
BusName=org.chromium.cras
ExecStartPre=/usr/bin/chmod -R 777 /var/run/cras
ExecStart=/usr/bin/sg audio /usr/bin/cras

[Install]
WantedBy=sound.target
