# krb5-appl-svn/krb5-eklogind@.service
[Unit]
Description=Encrypted Kerberized Remote Login service

[Service]
ExecStart=/usr/sbin/krb5-klogind -c -e
StandardInput=socket
