# krb5-appl-svn/krb5-telnetd@.service
[Unit]
Description=Kerberized Telnet server

[Service]
ExecStart=/usr/sbin/krb5-telnetd -h
StandardInput=socket
