# krb5-appl-svn/krb5-ftpd@.service
[Unit]
Description=Kerberized FTP Service

[Service]
ExecStart=/usr/sbin/krb5-ftpd -a
StandardInput=socket
