# krb5-appl-svn/krb5-kshd@.service
[Unit]
Description=Kerberized Remote Shell server

[Service]
ExecStart=/usr/sbin/krb5-kshd -5c
StandardInput=socket
