# krb5-appl-svn/krb5-telnetd.socket
[Unit]
Description=Kerberized Telnet socket

[Socket]
ListenStream=23
Accept=yes

[Install]
WantedBy=sockets.target
