# krb5-appl-svn/krb5-klogind@.service
[Unit]
Description=Kerberized Remote Login service

[Service]
ExecStart=/usr/sbin/krb5-klogind -c
StandardInput=socket
