# obelisk-devel-git/obelisk.service
[Unit]
Description=Obelisk server
Requires=sx-initchain.service
After=sx-initchain.service
OnFailureJobMode=ignore-requirements

[Service]
Type=simple
User=obelisk
Group=obelisk
ExecStart=/usr/bin/obelisk-server

[Install]
WantedBy=multi-user.target
