# adept-runtime-xilinx/20-digilent-usb.rules
# Digilent USB JTAG
ACTION=="add", SUBSYSTEMS=="usb", ATTRS{idVendor}=="0403", ATTRS{idProduct}=="6014", GROUP="uucp", MODE="0660" RUN+="/usr/bin/usbdrv -B %s{busnum} -A %s{devnum} -d"
ACTION=="add", SUBSYSTEMS=="usb", ATTRS{idVendor}=="0403", ATTRS{idProduct}=="6010", GROUP="uucp", MODE="0660" RUN+="/usr/bin/usbdrv -B %s{busnum} -A %s{devnum} -d"
