# systemd-mysqld-units/mysqld.service
[Unit]
Description=MySQL Server
After=network.target

[Service]
ExecStart=/usr/bin/mysqld --user=mysql --pid-file=/run/mysqld/mysqld.pid
ExecStop=/bin/kill -15 $MAINPID
PIDFile=/run/mysqld/mysqld.pid
Restart=always

[Install]
WantedBy=multi-user.target

