# miniupnp-git/minissdpd.service
[Unit]
Description=MiniSSDPd device discovery caching daemon
After=netctl.service

[Service]
Restart=always
EnvironmentFile=/etc/conf.d/minissdpd
ExecStart=/sbin/minissdpd -f -i $INTERFACE
PIDFile=/var/run/minissdpd.pid

[Install]
WantedBy=network.target
