# arielle-git/arielle.service
[Unit]
Description=Arielle IRC robot
Wants=network.target

[Service]
User=%I
EnvironmentFile=/etc/conf.d/arielle
ExecStart=/usr/bin/arielle $ARIELLE_ARGS

[Install]
WantedBy=multi-user.target
