# pommed-ng/gpomme.service
[Unit]
Description=GPomme Apple Hotkeys Daemon GUI

[Service]
User=root
ExecStart=/usr/bin/gpomme

[Install]
WantedBy=multi-user.target

