# opennebula/oned
#!/bin/sh -e
# chkconfig: 2345 20 80
# description: Starts the OpenNebula service
case "$1" in
        start)
                su - oneadmin /usr/bin/one start
                ;;
        stop)
                su - oneadmin /usr/bin/one stop
                ;;
        restart)
                su - oneadmin /usr/bin/one stop
                su - oneadmin /usr/bin/one start
                ;;
        *)
                echo "Usage: one {start|stop|restart}" >&2
                exit 1
                ;;

esac
exit 0
