# ipxripd/ipxd.service
[Unit]
Description=IPX RIP/SAP daemon

[Service]
Type=forking
ExecStart=/usr/sbin/ipxd
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
