# eibd/eibd.service
[Unit]
Description=EIB Daemon
After=network.target

[Service]
EnvironmentFile=/etc/conf.d/eibd.conf
ExecStart=/usr/bin/eibd -d -p /run/eibd.pid $EIBD_OPTS
Type=forking
PIDFile=/run/eibd.pid

[Install]
WantedBy=multi-user.target
