# spectrum/spectrum@.service
[Unit]
Description=Spectrum XMPP transport for %i
After=network.service prosody.service ejabberd.service

[Service]
Type=forking
PIDFile=/run/spectrum/%i.pid
ExecStart=/usr/bin/spectrumctl --su=spectrum -c /etc/spectrum/%i.cfg start
ExecStop=/usr/binspectrumctl --su=spectrum -c /etc/spectrum/%i.cfg stop
ExecReload=/usr/binspectrumctl --su=spectrum -c /etc/spectrum/%i.cfg reload

[Install]
WantedBy=multi-user.target

