# pvfs/pvfs2-client
#!/bin/sh
# Archlinux rc system startup file by andrew@lightningtoads.com
#
# description: pvfs2-client is the client component of PVFS2
#
# chkconfig: 345 35 55

# Source function library.
. /etc/rc.conf
. /etc/rc.d/functions

# executable
PVFS2CLIENT=/usr/sbin/pvfs2-client

# verify presence of client binary
if ! [ -x ${PVFS2CLIENT} ]; then
	echo "Error: could not find executable ${PVFS2CLIENT}"
	exit 1
fi

# See how we were called.
PID=`pidof -o %PPID /usr/sbin/pvfs2-client`
case "$1" in
  start)
	stat_busy "Starting PVFS2 client: "
	[ -z "$PID" ] && ${PVFS2CLIENT}
	if [ $? -gt 0 ]; then
		stat_fail
	else
		PID=`pidof -o %PPID /usr/sbin/pvfs2-client`
		echo $PID > /var/run/pvfs2-client.pid
		add_daemon pvfs2-client
		stat_done
	fi
	;;
  stop)
	stat_busy "Stopping PVFS2 client: "
	[ ! -z "$PID" ]  && kill $PID &> /dev/null
	if [ $? -gt 0 ]; then
		stat_fail
	else
		rm -f /var/run/pvfs2-client.pid
		rm_daemon pvfs2-client
		stat_done
	fi
	;;
  status)
        status pvfs2-client
	;;
  restart)
	$0 stop
	# give server time to die cleanly
	sleep 2
	$0 start
	;;
  *)
	echo "Usage: $0 {start|stop|status|restart}"
	exit 1
esac

exit 0

