# opensand-daemon/opensand-daemon.service
[Unit]
Description="Daemon for OpenSAND entity (sat, gw, st or ws)"
After=syslog.target network.target

[Service]
User=opensand
ExecStart=sand_daemon -c /etc/opensand/daemon.conf

[Install]
WantedBy=multi-user.target

