# microchip-mplab_bin/z010_mchp_tools.rules
# BxS 2011.11.18 Changed hotplugscript path
# BxS 2011.11.10 Fixed for udev 175
# PRJ 2010.01.26 Add reference to "usb" for Ubuntu.
# PRJ 2010.01.22 Attempt to further simplify rules files requirements.
# PRJ 2009.08.18 Rules file simplified.
# PRJ 2009.07.15 Rules file created.

ENV{hotplugscript}="/opt/microchip/mplabx/mplab_ide/mplablibs/modules/lib/mchplinusbdevice"

ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="04d8", MODE="0666", RUN+="$env{hotplugscript} add"

ACTION=="remove", SUBSYSTEM=="usb", RUN+="$env{hotplugscript} remove $env{PRODUCT}"
