# arpwatch-gentoo/arpwatch.service
[Unit]
Description=Network monitoring tool
After=network.target

[Service]
Type=forking
PIDFile=/run/arpwatch.pid
EnvironmentFile=/etc/conf.d/arpwatch
ExecStart=/usr/sbin/arpwatch $ARPWATCH_ARGS
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
