# logrotate-systemd/logrotate.service
[Unit]
Description=Rotate Logs

[Service]
Nice=19
IOSchedulingClass=2
IOSchedulingPriority=7
ExecStart=/usr/bin/logrotate /etc/logrotate.conf

[Install]
WantedBy=cron-daily.target
