# zarafa-licensed/zarafa-licensed.service
[Unit]
Description=Zarafa Collaboration Platform's License Server
After=network.target
Before=/usr/bin/zarafa-server

[Service]
Type=forking
GuessMainPID=no
Environment=LC_ALL=C LANG=C LD_LIBRARY_PATH=/opt/ubuntu1204libs/
ExecStart=/usr/bin/zarafa-licensed 2>/dev/null
TimeoutSec=65

[Install]
WantedBy=multi-user.target
