# chaosvpn/chaosvpn.service
[Unit]
Description=ChaosVPN

[Service]
ExecStart=/usr/sbin/chaosvpn
KillMode=process
Restart=always

[Install]
WantedBy=multi-user.target
