# logstash/logstash-web.service
[Unit]
Description=Logstash Web Server

[Service]
User=logstash
Group=logstash
ExecStart=/usr/bin/logstash web
PrivateTmp=true

[Install]
WantedBy=multi-user.target
