# ganglia/gmetad.service
[Unit]
Description=Ganglia Meta Daemon
After=network.target

[Service]
Type=forking
PIDFile=/run/gmetad.pid
ExecStart=/usr/sbin/gmetad -p /run/gmetad.pid

[Install]
WantedBy=multi-user.target
