# ganglia/gmond.service
[Unit]
Description=Ganglia Monitor Daemon
After=network.target

[Service]
ExecStart=/usr/sbin/gmond -f -p /run/gmond.pid

[Install]
WantedBy=multi-user.target
