# shadowsocks-libev-git/shadowsocks@.service
[Unit]
Description=Shadowsocks server
After=network.target

[Service]
Type=simple
PIDFile=/run/shadowsocks/%i.pid
ExecStart=/usr/bin/ss-server -f /run/shadowsocks/%i.pid -c /etc/shadowsocks/%i.json

[Install]
WantedBy=multi-user.target
