# shadowsocks-libev-git/shadowsocks.service
[Unit]
Description=Shadowsocks server
After=network.target

[Service]
Type=simple
PIDFile=/run/shadowsocks.pid
ExecStart=/usr/bin/ss-server -f /run/shadowsocks.pid -c /etc/shadowsocks/config.json

[Install]
WantedBy=multi-user.target
