# okimfpsdrv/okiscand.service
[Unit]
Description=OKI Scan Service

[Service]
ExecStart=/etc/init.d/okiscand start
ExecStop=/etc/init.d/okiscand stop
PIDFile=/var/run/okiscand.pid
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target

