# mit-zephyr/zhm.service
[Unit]
Description=Zephyr Host Manager
After=local-fs.target network.target

[Service]
RemainAfterExit=true
ExecStart=/usr/bin/zhm -d
ExecStop=/usr/bin/killall zhm
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=on-failure

[Install]
WantedBy=multi-user.target