# sympa/sympa-archived.service
[Unit]
Description=Web archives building for Sympa
BindTo=sympa.service
 
[Service]
User=sympa
Group=sympa
Type=forking
PIDFile=/var/lib/sympa/archived.pid
ExecStart=/usr/lib/sympa/archived.pl
