# sympa/sympa-task_manager.service
[Unit]
Description=Scheduled tasks processing for Sympa
BindTo=sympa.service
 
[Service]
User=sympa
Group=sympa
Type=forking
PIDFile=/srv/sympa/task_manager.pid
ExecStart=/opt/sympa/bin/task_manager.pl
