# sympa/sympa-archived.service
[Unit]
Description=Web archives building for Sympa
BindTo=sympa.service
 
[Service]
User=sympa
Group=sympa
Type=forking
PIDFile=/srv/sympa/archived.pid
ExecStart=/opt/sympa/bin/archived.pl
