# sympa/sympa-bounced.service
[Unit]
Description=Non-delivery reports analysis and storage for Sympa
BindTo=sympa.service
 
[Service]
User=sympa
Group=sympa
Type=forking
PIDFile=/var/lib/sympa/bounced.pid
ExecStart=/usr/lib/sympa/bounced.pl
