# sympa/sympa-bounced.service
[Unit]
Description=Non-delivery reports analysis and storage for Sympa
BindTo=sympa.service
 
[Service]
User=sympa
Group=sympa
Type=forking
PIDFile=/srv/sympa/bounced.pid
ExecStart=/opt/sympa/bin/bounced.pl
