# sympa/sympa-bulk.service
[Unit]
Description=Mail sending for Sympa
BindTo=sympa.service
 
[Service]
User=sympa
Group=sympa
Type=forking
PIDFile=/srv/sympa/bulk.pid
ExecStart=/opt/sympa/bin/bulk.pl
