# sympa/sympa.service
[Unit]
Description=Main Sympa daemon
Wants=sympa-bulk.service
Wants=sympa-archived.service
Wants=sympa-bounced.service
Wants=sympa-task_manager.service
 
[Service]
User=sympa
Group=sympa
Type=forking
PIDFile=/var/lib/sympa/sympa.pid
ExecStart=/usr/lib/sympa/sympa.pl
 
[Install]
WantedBy=multi-user.target
