# redeclipse-server/redeclipse.service
[Unit]
Description=Redeclipse Server
After=network.target

[Service]
User=nobody
Type=simple
Restart=on-abort
StandardOutput=syslog
ExecStart=/usr/bin/redeclipse-server -h/etc/redeclipse

[Install]
WantedBy=multi-user.target
