# cassandra/cassandra.service
[Unit]
Description=Cassandra NoSQL database
After=network.target

[Service]
Type=forking
User=cassandra
LimitMEMLOCK=infinity
LimitNOFILE=100000
LimitNPROC=32768
LimitAS=infinity
PIDFile=/var/run/cassandra/cassandra.pid
ExecStart=/usr/bin/cassandra -p /var/run/cassandra/cassandra.pid
WorkingDirectory=/var/lib/cassandra

[Install]
WantedBy=multi-user.target
