# pgl/pgl-update.service
[Unit]
Description=Update PeerGuardian Linux blocklists
Wants=network.target

[Service]
EnvironmentFile=/usr/lib/pgl/pglcmd.main
ExecStart=$CMD_PATHNAME update

[Install]
WantedBy=multi-user.target