# pgl/pgl.service
[Unit]
Description=PeerGuardian Linux - an IP Blocker
Wants=network.target

[Service]
Type=oneshot
EnvironmentFile=/usr/lib/pgl/pglcmd.main
ExecStart=$CMD_PATHNAME start
ExecStop=$CMD_PATHNAME stop_quick
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target