# radeon-modeset-service/radeon-modeset.service
[Unit]
Description=Enable Radeon modeset
Before=display-manager.service

[Service]
Type=oneshot
ExecStart=/sbin/modprobe -r radeon drm ; /sbin/modprobe radeon modeset=1

[Install]
WantedBy=multi-user.target
