# pulseaudio-systemd/pulseaudio.service
[Unit]
Description=PulseAudio Sound Server in system-wide mode

[Service]
Type=forking
PIDFile=/var/run/pulse/pid
ExecStart=/usr/bin/pulseaudio --system --disallow-module-loading=1 \
	--disallow-exit=1 --disable-shm=1 --fail=1 --daemonize

[Install]
WantedBy=multi-user.target
