# systemd-dropbox-daemon/dropbox@.service
[Unit]
Description=Dropbox

[Service]
Type=simple
#EnvironmentFile=/etc/conf.d/dropboxd
ExecStart=/usr/bin/dropboxd
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=always
User=%I

[Install]
WantedBy=multi-user.target
