# systemd-dropbox-daemon/dropboxd.service
[Unit]
Description=Dropbox Daemon as a user

[Service]
Type=oneshot
ExecStart=/etc/rc.d/dropboxd start
ExecStop=/etc/rc.d/dropboxd stop
ExecRestart=/etc/rc.d/dropboxd restart
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target

