# qubes-vm-gui/qubes-gui-agent.service
[Unit]
Description=Qubes GUI Agent
After=qubes-meminfo-writer.service qubes-dvm.service qubes-misc-post.service

[Service]
# pretend that user is at local console
ExecStartPre=/usr/bin/touch /var/run/console/user
# clean env
ExecStartPre=/usr/bin/rm -f /tmp/qubes-session-env /tmp/qubes-session-waiter
ExecStart=/usr/bin/qubes_gui $GUI_OPTS
StandardOutput=syslog
Environment=DISPLAY=:0
EnvironmentFile=-/var/run/qubes-service-environment

[Install]
WantedBy=multi-user.target
