# asterisk18/asterisk.service
[Unit]
Description=A complete PBX solution
After=local-fs.target

[Service]
Type=forking
ExecStart=/usr/sbin/asterisk -G asterisk -U asterisk
ExecStop=kill /var/run/asterisk/asterisk.pid &> /dev/null

[Install]
WantedBy=multi-user.target
