# syncthing/syncthing@.service
[Unit]
Description=Syncthing service for %i
After=network.target

[Service]
User=%i
ExecStart=/usr/bin/syncthing -home=/home/%i/.syncthing
Restart=on-abort

[Install]
WantedBy=multi-user.target
