# cpuminer-git/cpuminer.service
[Unit]
Description=Multi-threaded CPU miner for Bitcoin and Litecoin
After=network.target

[Service]
Type=simple
PIDFile=/run/minerd.pid
ExecStart=/usr/bin/minerd --config /etc/cpuminer.json
ExecStop=/bin/kill -TERM $MAINPID
CPUShares=1

[Install]
WantedBy=multi-user.target