# obdi-git/obdi-worker.service
[Unit]
Description=Obdi

[Service]
Type=simple
Restart=always
PIDFile=/var/run/obdi-worker.pid
ExecStart=/usr/bin/obdi-worker

[Install]
WantedBy=multi-user.target

