# obdi-git/obdi.service
[Unit]
Description=Obdi

[Service]
Type=simple
Restart=always
PIDFile=/var/run/obdi.pid
ExecStart=/usr/bin/obdi

[Install]
WantedBy=multi-user.target

