# buddycloud-http-api-git/buddycloud-http-api.service
[Unit]
Description=buddycloud HTTP API (git version)
After=network.target

[Service]
User=nobody
WorkingDirectory=/opt/buddycloud-http-api
Type=simple
ExecStart=/bin/bash -c '/usr/bin/node -- ./server.js >>/var/log/buddycloud-http-api.log 2>&1'

[Install]
WantedBy=multi-user.target
