# kafka/kafka.service
[Unit]
Description=Highly reliable distributed coordination server
After=network.target

[Service]
Type=simple
User=root
SyslogIdentifier=kafka
Environment="KAFKA_LOG4J_OPTS=-Dlog4j.configuration=file:/etc/kafka/log4j.properties"
ExecStart=/opt/kafka_2.10-0.8.1.1/bin/kafka-run-class.sh -name kafkaServer -loggc kafka.Kafka /etc/kafka/server.properties
ExecStop=/opt/kafka_2.10-0.8.1.1/bin/kafka-server-stop.sh

[Install]
WantedBy=multi-user.target