# kafka/systemd_kafka.service
[Unit]
Description=A high-throughput distributed messaging system
Requires=kafka-zookeeper.service
After=network.target kafka-zookeeper.service

[Service]
Type=forking
SuccessExitStatus=143
ExecStart=/usr/bin/kafka-server-start.sh -daemon /etc/kafka/server.properties
ExecStop=/usr/bin/kafka-server-stop.sh

[Install]
WantedBy=multi-user.target
