# kafka/systemd_kafka-zookeeper.service
[Unit]
Description=An open-source server which enables highly reliable distributed coordination
After=network.target
StopWhenUnneeded=true

[Service]
Type=forking
SuccessExitStatus=143
ExecStart=/usr/bin/zookeeper-server-start.sh -daemon /etc/kafka/zookeeper.properties
ExecStop=/usr/bin/zookeeper-server-stop.sh

[Install]
WantedBy=multi-user.target
