# vmware-workstation/vmware-usb.service
[Unit]
Description=VMware USB Arbitration Service
After=vmware-authentication.service
PartOf=vmware.target

[Service]
Type=simple
RemainAfterExit=yes
ExecStart=/opt/vmware/bin/vmware-usbarbitrator
ExecStop=/opt/vmware/bin/vmware-usbarbitrator --kill

[Install]
WantedBy=vmware.target
