# vmware-workstation/vmware-vmblock.service
[Unit]
Description=VMware Blocking File System
After=vmware-vmci.service
PartOf=vmware.target

[Service]
Type=simple
RemainAfterExit=yes
ExecStart=/sbin/modprobe vmblock
ExecStop=/sbin/modprobe -r vmblock

[Install]
WantedBy=vmware.target
