# vmware-workstation/vmware-authentication.service
[Unit]
Description=VMware Authentication Daemon
After=vmware-vmnet.service
PartOf=vmware.target

[Service]
Type=simple
RemainAfterExit=yes
ExecStart=/opt/vmware/sbin/vmware-authdlauncher

[Install]
WantedBy=vmware.target
