# vmware-workstation/vmware.target
[Unit]
Description=VMware Workstation Service(s)
Requires=vmware-vmmon.service
Wants=vmware-vmci.service
Wants=vmware-vmsock.service
Wants=vmware-vmblock.service
Wants=vmware-vmnet.service
Wants=vmware-authentication.service
Wants=vmware-usb.service

[Install]
WantedBy=multi-user.target
