# vmware-workstation/vmware-vmmon.service
[Unit]
Description=VMware Virtual Machine Monitor
PartOf=vmware.target

[Service]
Type=simple
RemainAfterExit=yes
ExecStart=/sbin/modprobe vmmon
ExecStop=/sbin/modprobe -r vmmon

[Install]
WantedBy=vmware.target
