# vmware-workstation/vmware-vmci.service
[Unit]
Description=VMware Virtual Machine Communication Interface
After=vmware-vmmon.service
PartOf=vmware.target

[Service]
Type=simple
RemainAfterExit=yes
ExecStart=/sbin/modprobe vmci
ExecStop=/sbin/modprobe -r vmci

[Install]
WantedBy=vmware.target
