# vmware-workstation/vmware-vmsock.service
[Unit]
Description=VMware Virtual Machine Communication Interface Socket Family
Requires=vmware-vmci.service
After=vmware-vmci.service
PartOf=vmware.target

[Service]
Type=simple
RemainAfterExit=yes
ExecStart=/sbin/modprobe vsock
ExecStop=/sbin/modprobe -r vsock

[Install]
WantedBy=vmware.target
