# remctl/remctl@.service
[Unit]
Description=remctl command execution service
After=syslog.target

[Service]
ExecStart=/usr/sbin/remctld
StandardInput=socket
