# remctl/remctl@.service
[Unit]
Description=remctl command execution service
Documentation=man:remctld(8)

[Service]
ExecStart=/usr/sbin/remctld
StandardInput=socket
