# remctl/remctl@.service
[Unit]
Description=remctl command execution service

[Service]
ExecStart=/usr/sbin/remctld
StandardInput=socket
