# remctl/remctl.socket
[Unit]
Description=remctl command execution service
Conflicts=remctl.service

[Socket]
ListenStream=4373
Accept=yes

[Install]
WantedBy=sockets.target
