# lsi-msm/lsi_mrdsnmpd.service
[Unit]
Description=LSI MegaRAID Storage Manager SMNP Agent
PartOf=lsi_msm.service
BindsTo=lsi_msm.service
Requires=snmpd.service

[Service]
Type=forking
ExecStart=/usr/bin/lsi_mrdsnmpagent -c /etc/snmp/snmpd.conf

[Install]
WantedBy=multi-user.target
