# linux-gpib/gpib.rules
# linux-gpib udev rules

# NI GPIB USB

# NI GPIB USB B
SUBSYSTEM=="usb", ACTION=="add", ATTR{idVendor}=="3923", ATTR{idProduct}=="702b", RUN+="/sbin/fxload -D %N -t fx -I /usr/share/gpib-firmware/ni_gpib_usb_b/niusbb_firmware.hex -s /usr/share/gpib-firmware/ni_gpib_usb_b/niusbb_loader.hex"
# Note that 702b becomes 702a after firmware is loaded
SUBSYSTEM=="usb", ACTION=="add", ATTR{idVendor}=="3923", ATTR{idProduct}=="702a", RUN+="/usr/sbin/gpib_config"

# NI GPIB USB HS
SUBSYSTEM=="usb", ACTION=="add", ATTR{idVendor}=="3923", ATTR{idProduct}=="709b", RUN+="/usr/sbin/gpib_config"

# Agilent GPIB USB

# Agilent 82357a
SUBSYSTEM=="usb", ACTION=="add", ATTR{idVendor}=="957", ATTR{idProduct}=="7", RUN+="/sbin/fxload -D %N -I /usr/share/gpib-firmware/agilent_82357a/82357a_fw.hex
# Note that 7 becomes 107 after firmware is loaded
SUBSYSTEM=="usb", ACTION=="add", ATTR{idVendor}=="957", ATTR{idProduct}=="107", RUN+="/usr/sbin/gpib_config"

# Agilent 82357b
SUBSYSTEM=="usb", ACTION=="add", ATTR{idVendor}=="957", ATTR{idProduct}=="7", RUN+="/sbin/fxload -D %N -t fx2 -I /usr/share/gpib-firmware/agilent_82357a/measat_releaseX1.8.hex
# Note that 7 becomes 107 after firmware is loaded
SUBSYSTEM=="usb", ACTION=="add", ATTR{idVendor}=="957", ATTR{idProduct}=="107", RUN+="/usr/sbin/gpib_config"

# Devices
KERNEL=="gpib/*",       MODE="0660", GROUP="gpib"
KERNEL=="gpib[0-9]*",   MODE="0660", GROUP="gpib"

