# actkbd/actkbd.service
[Unit]
Description=Actkbd: Daemon for X-independent shortcuts
Wants=setkeycodes.service
Requires=setkeycodes.service

[Service]
Type=forking
ExecStart=/usr/bin/actkbd -D
RestartSec=1
User=root

[Install]
WantedBy=multi-user.target
