# groestlcoin-daemon/groestlcoind@.service
[Unit]
Description=GroestlCoin Daemon
After=network.target 

[Service]
User=%i
Group=users
Type=simple
ExecStart=/usr/bin/groestlcoind

[Install]
WantedBy=multi-user.target

