# bitwrk-git/bitwrk.service
[Unit]
Description=Bitwrk Client Daemon

[Service]
EnvironmentFile=/etc/conf.d/bitwrk.conf
ExecStart=/usr/bin/nice /usr/bin/bitwrk-client $OPTIONS 2>&1 &
ExecStop=/bin/pkill /usr/bin/bitwrk-client

[Install]
WantedBy=multi-user.target
