# powertop-autotune-systemd/powertop-autotune.service
# Credits: https://wiki.archlinux.org/index.php/Powertop#Tips_and_tricks
[Unit]
Description=Powertop autotune

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/bin/powertop --auto-tune
Environment="TERM=xterm"

[Install]
WantedBy=multi-user.target
