# radsecproxy/radsecproxy.service
[Unit]
Description=generic RADIUS proxy
After=network.target

[Service]
Type=forking
Restart=always
ExecStartPre=/usr/bin/radsecproxy -p
ExecStart=/usr/bin/radsecproxy

[Install]
WantedBy=multi-user.target
