# echovnc-server/echoserver.service
[Unit]
Description=echoserver EchoVNC Realy
After=network.target

[Service]
ExecStart=/usr/bin/echoserver -c /etc/echoserver.conf
Type=forking

[Install]
WantedBy=multi-user.target
