# influxdb/influxdb.service
[Unit]
Description=InfluxDB

[Service]
User=influxdb
Group=influxdb
Type=install
ExecStart=/usr/bin/influxdb -config=/etc/influxdb.conf -pidfile=/run/influxdb/influxdb.pid

[Install]
WantedBy=multi-user.target
