# ventrilo-server/ventrilo-server@.service
[Unit]
Description=Ventrilo Server
After=network.target

[Service]
Type=forking
User=%i
PIDFile=/opt/ventsrv/ventrilo_srv.pid
ExecStart=/opt/ventsrv/ventrilo_srv -d
WorkingDirectory=/opt/ventsrv

[Install]
WantedBy=multi-user.target

