# iplist/ipblock.service
[Unit]
Description=a Linux application for blocking connections to and from a specified range of hosts using the netfilter netlink-queue library.
After=ufw.service

[Service]
Type=forking
PIDFile=/var/run/iplist.pid
ExecStart=/usr/sbin/ipblock -s
ExecStop=/usr/sbin/ipblock -d

[Install]
WantedBy=multi-user.target
