# macfanctld-git/macfanctld.service
[Unit]
Description=Macfanctld

[Service]
Type=simple
ExecStart=/usr/bin/macfanctld
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=always

[Install]
WantedBy=multi-user.target
