# thinkfan/thinkfan.service
[Unit]
Description=Thinkfan Daemon - simple and lightweight fan control program
After=basic.target

[Service]
ExecStart=/usr/sbin/thinkfan -n -s5

[Install]
WantedBy=multi-user.target
