# thinkfan/thinkfan.service
[Unit]
Description=Thinkfan Daemon - simple and lightweight fan control program
After=basic.target

[Service]
ExecStart=/usr/bin/thinkfan -n -s5 -q

[Install]
WantedBy=multi-user.target
