# emulationstation-git-unstable-rpi/emulationstation.service
[Unit]
Description = Starts an instance of Emulationstation
After = remote-fs.target

[Service]
User = emulationstation
Group = emulationstation
Type = simple
ExecStart = /usr/bin/emulationstation
Restart = on-failure

[Install]
WantedBy = multi-user.target

