# cloud9/cloud9.service
[Unit]
Description=Cloud9 server
After=syslog.target network.target

[Service]
Type=simple
User=cloud9
Group=cloud9
WorkingDirectory=/var/lib/cloud9
ExecStart=/usr/share/webapps/cloud9/bin/cloud9.sh
Restart=always

[Install]
WantedBy=multi-user.target
