# phddns/phddns.service
[Unit]
Description=
After=network.target

[Service]
Type=forking
PIDFile=/var/run/phddns.pid
ExecStartPre=/bin/rm -f /var/run/phddns.pid
ExecStart=/usr/bin/phddns --daemon --config /etc/phddns.conf
Restart=on-abort

[Install]
WantedBy=multi-user.target
