# mod_gearman/mod-gearman-worker.service
[Unit]
Description=Mod Gearman worker

[Service]
User=gearman
ExecStart=/usr/bin/mod_gearman_worker --config=/etc/mod_gearman/mod_gearman_worker.conf

[Install]
WantedBy=multi-user.target
