# zabbix-server-mysql/zabbix-agentd.service
[Unit]
Description=Zabbix Agent
After=syslog.target network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/zabbix_agentd
RemainAfterExit=yes
PIDFile=/var/run/zabbix/zabbix_agentd.pid

[Install]
WantedBy=multi-user.target
