# xen-4.3/xendomU.service
[Unit]
Description=xen domU service for domain %i
Requires=xenstored.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStartPre=/usr/bin/xenstore-write "/local/domain/0/name" "Domain-0
ExecStart=/usr/sbin/xl create /etc/xen/conf/%i.cfg
ExecStop=/usr/sbin/xl shutdown %i

[Install]
WantedBy=multi-user.target
