# gitlab-ci/gitlab-ci-unicorn.service
[Unit]
Description=GitLab CI Unicorn Server

[Service]
User=gitlab_ci
WorkingDirectory=/usr/share/webapps/gitlab-ci
Environment=RAILS_ENV=production
SyslogIdentifier=gitlab-ci-unicorn
PIDFile=/var/run/gitlab-ci/unicorn.pid

ExecStart=/usr/bin/bundle exec "unicorn_rails -c /usr/share/webapps/gitlab-ci/config/unicorn.rb -E production"
ExecStop=/usr/bin/kill -QUIT $MAINPID
ExecReload=/usr/bin/kill -USR2 $MAINPID

[Install]
WantedBy=gitlab-ci.target
