# kippo/kippo.service
[Unit]
Description=Kippo Honey Pot
After=network.target

[Service]
Type=forking
WorkingDirectory=/opt/kippo
ExecStart=/opt/kippo/start.sh

[Install]
WantedBy=multi-user.target
