# sonarr-git/sonarr.service
[Unit]
Description=Sonarr Service
After=network.target

[Service]
User=sonarr
Group=sonarr
Restart=always
RestartSec=5
Type=simple
ExecStart=/opt/sonarr -data=/var/lib/sonarr
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
