# gogs/gogs.service
[Unit]
Description=Gogs(Go Git Service)
After=network.target

[Service]
Type=simple
ExecStart=/var/lib/gogs/gogs web 

[Install]
WantedBy=multi-user.target
