# seafile-server/seafile-server@.service
[Unit]
Description=Seafile is an open-source Dropbox replacement.
After=syslog.target network.target

[Service]
Type=forking
User=seafile
WorkingDirectory=%h/%i
ExecStart=/usr/bin/seafile-admin start --fastcgi
ExecStop=/usr/bin/seafile-admin stop
PIDFile=%h/%i/seafile-data/pids/seaf-server.pid

[Install]
WantedBy=multi-user.target
