# git-git/git-daemon.service
[Unit]
Description=Git Daemon Instance

[Service]
User=git
ExecStart=-/usr/lib/git-core/git-daemon --export-all
