# ibconnect-git/ibconnect.service
[Unit]
Description=IB Connect

[Service]
Type=simple
User=root
Group=root

TimeoutStopSec=6
KillMode=control-group
KillSignal=SIGTERM
SendSIGHUP=no
SendSIGKILL=yes

EnvironmentFile=/etc/ibconnect
ExecStart=/usr/bin/ibcd

[Install]
WantedBy=multi-user.target
