# init-headphone/init-headphone.service
[Unit]
Description=Reactivate headphones after sleep
After=sleep.target

[Service]
Type=oneshot
ExecStart=/usr/sbin/init-headphone

[Install]
WantedBy=sleep.target

