# granola-connect/granola-connect.service
[Unit]
Description=Granola Connection Daemon
After=granola.service

[Service]
ExecStart=/usr/sbin/granola-connect
