# ndnlp-git/ndnld.service
[Unit]
Description=Named Data Networking Link Protocol
After=network.target
Requires=ccnd.service

[Service]
#Type=forking
ExecStart=/usr/bin/ndnld

[Install]
WantedBy=multi-user.target
