# systemd-memcached-units/memcached.service
[Unit]
Description=Memcached Server
After=network.target

[Service]
EnvironmentFile=/etc/conf.d/memcached
ExecStart=/usr/bin/memcached -u $MEMCACHED_USER $MEMCACHED_ARGS

[Install]
WantedBy=multi-user.target
