# xpra-systemd/xpra-client@.service
[Unit]
Description=xpra client displaying :%i
Requires=xpra@%i.service
After=xpra@%i.service

[Service]
Type=simple
Restart=on-failure
# The server needs some time to start up. Give it chances.
RestartSec=1s
StartLimitInterval=10
StartLimitBurst=9
Environment="DISPLAY=:0"
ExecStart=/usr/bin/xpra attach :%i

[Install]
WantedBy=default.target
