# oscam/oscam.service
[Unit]
Description=An Open Source Conditional Access Module software

[Service]
Type=forking
ExecStart=/usr/bin/oscam -b -c /etc/oscam

[Install]
WantedBy=multi-user.target
