# plymouth-lite/plymouth-lite-start.service
[Unit]
Description=Show Plymouth-lite Start Screen
DefaultDependencies=no
After=systemd-vconsole-setup.service
Before=sysinit.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/ply-image /usr/share/plymouth/splash.png > /dev/null 2>&1
