# plymouth-lite/plymouth-lite-poweroff.service
[Unit]
Description=Show Plymouth-lite Poweroff Screen
DefaultDependencies=no
Before=systemd-poweroff.service
After=plymouth-lite-start.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/ply-image /usr/share/plymouth/poweroff.png > /dev/null 2>&1
