# plymouth-lite/plymouth-lite-halt.service
[Unit]
Description=Show Plymouth-lite Halt Screen
DefaultDependencies=no
Before=systemd-halt.service
After=plymouth-lite-start.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/ply-image /usr/share/plymouth/halt.png > /dev/null 2>&1
