# plymouth-lite/plymouth-lite-reboot.service
[Unit]
Description=Show Plymouth-lite Reboot Screen
DefaultDependencies=no
Before=systemd-reboot.service
After=plymouth-lite-start.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/ply-image /usr/share/plymouth/reboot.png > /dev/null 2>&1
