# spindown/spindown.service
[Unit]
Description=Spindown Daemon

[Service]
ExecStart=/usr/bin/spindownd -d -f /run/spindown/spindownd.fifo -c /etc/spindown.conf

[Install]
WantedBy=multi-user.target
