# samsung-tools/samsung-tools.service
[Unit]
Description=Tools for Samsung laptops

[Service]
ExecStart=/usr/bin/dbus-send --system --print-reply=literal --dest=org.voria.SamsungTools.System / org.voria.SamsungTools.System.SetInitialDevicesStatus
Type=forking

[Install]
WantedBy=multi-user.target
