# ttdnsd-git/ttdnsd.service
[Unit]
Description=Tor TCP DNS Daemon
Requires=tor.service
After=tor.service

[Service]
Type=forking
PIDFile=/run/ttdnsd.pid
EnvironmentFile=/etc/conf.d/ttdnsd
ExecStart=/usr/sbin/ttdnsd -P /run/ttdnsd.pid -f /etc/ttdnsd.conf -C /var/lib/ttdnsd $DEFAULTS
Restart=always

[Install]
WantedBy=multi-user.target
