# palm-novacom-lite/novacomd.service
[Unit]
Description=Novacom service for webOS development connectivity

[Service]
Type=forking
Environment=LD_LIBRARY_PATH=/opt/Palm/novacom
EnvironmentFile=/etc/conf.d/novacomd
GuessMainPID=yes
ExecStart=/usr/bin/novacomd $NOVACOMD_OPTIONS

[Install]
WantedBy=multi-user.target
