# plymouth/kdm-plymouth.service
[Unit]
Description=K Display Manager
After=systemd-user-sessions.service

Conflicts=plymouth-quit.service
After=plymouth-quit.service

[Service]
ExecStart=/usr/bin/kdm -nodaemon
StandardOutput=syslog

[Install]
Alias=display-manager.service
