# perspectives-server-git/perspectives.service
[Unit]
Description=Perspectives notary server
After=network.target

[Service]
Type=forking
PIDFile=/opt/perspectives/perspectives.pid
User=perspectives
ExecStart=usr/bin/python2 /opt/perspectives/notary_http.py
WorkingDirectory=/opt/perspectives

[Install]
WantedBy=multi-user.target

