# oscam-svn/oscam.service
[Unit]
Description=An Open Source Conditional Access Module software

[Service]
Type=forking
EnvironmentFile=/etc/conf.d/oscam
ExecStart=/usr/bin/oscam $OSCAM_ARGS

[Install]
WantedBy=multi-user.target
