# oscam-svn/oscam.service
[Unit]
Description=The Open Source Conditional Access Module daemon

[Service]
Type=forking
ExecStart=/usr/bin/oscam -b -c /etc/oscam

[Install]
WantedBy=multi-user.target
