# wview/wview.service
[Unit]
Description=WView Weather Station
After=network.target remote-fs.target

[Service]
Type=forking
PIDFile=/var/wview/wviewd.pid
ExecStart=/usr/bin/wviewctl start
ExecStop=/usr/bin/wviewctl stop
PrivateTmp=true

[Install]
WantedBy=multi-user.target
