# mcollective-2.2/mcollective.service
[Unit]
Description=The Marionette Collective
After=network.target

[Service]
Type=forking
StandardOutput=syslog
StandardError=syslog
ExecStart=/usr/bin/mcollectived --config=/etc/mcollective/%I.cfg --pidfile=/run/mcollective-%I.pid
ExecReload=/bin/kill -USR1 $MAINPID
PIDFile=/run/mcollective-%I.pid

[Install]
WantedBy=multi-user.target
