# jsession/jsession.service
[Unit]
Description=jsession X session manager

[Service]
Type=dbus
BusName=prog.jsession
ExecStart=/usr/bin/jsessiond
KillMode=process

[Install]
WantedBy=multi-user.target
