# openvswitch-git/ovs-vswitchd.service
[Unit]
Description    = Open vSwitch Daemon
Documentation  = man:ovs-vswitchd
Wants          = network.target
Before         = network.target
Requires       = ovsdb-server.service
After          = ovsdb-server.service

[Service]
Type=forking
ExecStartPre   = /sbin/modprobe openvswitch
ExecStart      = /usr/sbin/ovs-vswitchd \
                --pidfile=/run/openvswitch/ovs-vswitchd.pid \
                --detach


[Install]
WantedBy       = multi-user.target
