# openvswitch-git/ovsdb-server.service
[Unit]
Description    =  Open vSwitch Database Server
Documentation  =  man:ovsdb-server

[Service]
Type=forking
Type=forking
ExecStart=/usr/sbin/ovsdb-server \
             --remote=punix:/run/openvswitch/db.sock \
             --pidfile=/run/openvswitch/ovsdb-server.pid \
             --detach

ExecStartPost=/usr/bin/ovs-vsctl --no-wait init
