# devdocs-git/devdocs.service
[Unit]
Description=Devdocs instance
Requires=network.target
After=network.target

[Service]
User=devdocs
Group=devdocs
PrivateTmp=true
PrivateDevices=true
ProtectSystem=full
ProtectHome=true
WorkingDirectory=/usr/share/webapps/devdocs

ExecStart=/usr/bin/bundle exec rackup

[Install]
WantedBy=multi-user.target
