# libcec-daemon/libcec-daemon.service
[Unit]
Description=CEC Daemon using libCEC
After=syslog.target network.target

[Service]
Type=simple
ExecStart=/usr/bin/libcec-daemon
ExecStop=/usr/bin/killall -q -TERM libcec-daemon ; /usr/bin/killall -q -TERM libcec-daemon
KillMode=none

[Install]
WantedBy=multi-user.target
