# foundationdb/foundationdb.service
[Unit]
Description=Distributed key/value store with ACID transactions
After=network.target

[Service]
ExecStart=/usr/sbin/fdbmonitor --conffile /etc/foundationdb/foundationdb.conf --lockfile /var/run/fdbmonitor.pid --daemonize
PidFile=/var/run/fdbmonitor.pid
Type=forking

[Install]
WantedBy=multi-user.target
