# chinadns/chinadns.service
[Unit]
Description=ChinaDNS Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/chinadns

[Install]
WantedBy=multi-user.target
