# taskd/taskd.service
[Unit]
Description=A secure server providing multi-client access to task data
Documentation=man:taskd(1)
After=network.target

[Service]
EnvironmentFile=/etc/conf.d/taskd
ExecStart=/usr/bin/taskd server

[Install]
WantedBy=multi-user.target
