# taskd/taskd.service
[Unit]
Description=A secure server providing multi-client access to task data
Documentation=man:taskd(1)
After=network.target

[Service]
EnvironmentFile=/etc/conf.d/taskd
Type=forking
PIDFile=/var/lib/taskd/taskd.pid
ExecStart=/usr/bin/taskd server --daemon
# Ideally this is run as a non-root user, but at the very least we can do this hardening
ReadOnlyDirectories=/
ReadWriteDirectories=/var/lib/taskd

[Install]
WantedBy=multi-user.target
